# sigma/ai-agent-jadepuffer-autonomously-deploys-ransomware-in-novel-attack-SIGMA.yaml
title: Langflow RCE Attempt via Run API
id: 2a9b4c1d-8e7f-4a3b-9c1d-5e6f7a8b9c0d
status: experimental
description: Detects web requests attempting to exploit Langflow via the /api/v1/run/ endpoint, potentially indicating an RCE attempt similar to CVE-2025-3248.
references:
  - https://www.securityweek.com/agentic-ai-used-to-conduct-ransomware-attack-via-langflow/
author: cyber.netsecops.io
date: 2026-07-03
tags:
  - attack.t1190
level: high
logsource:
  category: webserver
detection:
  selection:
    cs-uri-query|contains: '/api/v1/run/'
  condition: selection
---
title: PostgreSQL Database Dump Execution
id: 3b8c5d2e-9f6e-4b2a-8d0c-4e5f6a7b8c9d
status: experimental
description: Detects the execution of the 'pg_dump' utility, which can be used by attackers to exfiltrate database contents.
references:
  - https://www.securityweek.com/agentic-ai-used-to-conduct-ransomware-attack-via-langflow/
author: cyber.netsecops.io
date: 2026-07-03
tags:
  - attack.t1213
level: medium
logsource:
  category: process_creation
  product: linux
detection:
  selection:
    Image|endswith: '/pg_dump'
  condition: selection
---
title: Cron Persistence File Modification
id: 546db805-acd3-44d2-aa73-9946dcb73195
status: experimental
description: Detects the creation or modification of files in cron directories, a common technique for establishing persistence on Linux systems.
references:
  - https://www.securityweek.com/agentic-ai-used-to-conduct-ransomware-attack-via-langflow/
author: cyber.netsecops.io
date: 2026-07-03
tags:
  - attack.t1053.003
level: medium
logsource:
  category: file_event
  product: linux
detection:
  selection:
    TargetFilename|startswith:
      - '/var/spool/cron/'
      - '/etc/cron.d/'
      - '/etc/cron.hourly/'
      - '/etc/cron.daily/'
      - '/etc/cron.weekly/'
      - '/etc/cron.monthly/'
  condition: selection