# sigma/china-linked-sprysocks-backdoor-evolves-with-new-windows-variants-SIGMA.yaml
title: SprySOCKS Scheduled Task Creation for Persistence
id: 8e1d9c2b-5a4f-4c6b-9d8e-e0f1a2b3c4d5
status: experimental
description: Detects the creation of a scheduled task, a persistence method used by SprySOCKS malware.
references:
  - https://thehackernews.com/2026/06/china-linked-sprysocks-backdoor-expands.html
author: cyber.netsecops.io
date: 2026-06-16
tags:
  - attack.t1053.005
  - attack.persistence
  - attack.execution
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
  condition: selection