# sigma/china-apt-mirrorface-expands-targeting-to-eu-diplomats-SIGMA.yaml
title: Suspicious LNK File Creation in Startup Folder
id: 9f8e7d6c-5b4a-4987-b6a5-4321fedcba98
status: experimental
description: Detects the creation of LNK files in user or system startup directories, a common persistence technique used by APTs.
references:
  - https://www.eset.com/int/about/newsroom/press-releases/research/eset-releases-latest-apt-report-china-aligned-groups-expand-targeting-iran-advances-diplomatic-esp/
author: cyber.netsecops.io
date: 2026-04-27
tags:
  - attack.t1547.009
  - attack.persistence
level: medium
logsource:
  category: file_event
  product: windows
detection:
  selection:
    TargetFilename|endswith: '.lnk'
    TargetFilename|contains:
      - '\AppData\Roaming\Microsoft\Windows\Start Menu\Programs\Startup\'
      - '\ProgramData\Microsoft\Windows\Start Menu\Programs\StartUp\'
  condition: selection
falsepositives:
  - Legitimate software installers creating startup shortcuts.
---
title: Suspicious Registry Run Key Creation
id: 1a2b3c4d-5e6f-7a8b-9c0d-1e2f3a4b5c6d
status: experimental
description: Detects modification of the registry run keys, a common persistence technique used by malware and APTs.
references:
  - https://www.eset.com/int/about/newsroom/press-releases/research/eset-releases-latest-apt-report-china-aligned-groups-expand-targeting-iran-advances-diplomatic-esp/
author: cyber.netsecops.io
date: 2026-04-27
tags:
  - attack.t1547.001
  - attack.persistence
level: medium
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    TargetObject|contains:
      - 'Software\Microsoft\Windows\CurrentVersion\Run'
      - 'Software\Microsoft\Windows\CurrentVersion\RunOnce'
    EventType: 'SetValue'
  condition: selection
falsepositives:
  - Legitimate software installers and updaters.