# sigma/jdownloader-website-hacked-in-supply-chain-attack-to-distribute-rat-SIGMA.yaml
title: JDownloader Linux Malware Persistence via profile.d
id: 8d7c6b5a-4f3e-4d2c-1b0a-987654321fed
status: experimental
description: Detects the creation of a malicious shell script in /etc/profile.d, a technique used by the Linux variant of the JDownloader malware for persistence.
references:
  - https://securityaffairs.co/wordpress/162985/malware/jdownloader-site-malware-supply-chain-attack.html
author: cyber.netsecops.io
date: 2026-05-10
tags:
  - attack.t1546.004
  - attack.persistence
level: high
logsource:
  category: file_event
  product: linux
detection:
  selection:
    TargetFilename: '/etc/profile.d/systemd.sh'
  condition: selection
---
title: JDownloader Linux Malware Masquerading as upowerd
id: 1e2f3a4b-5c6d-7e8f-9a0b-1c2d3e4f5a6b
status: experimental
description: 'Detects a process named upowerd running from a non-standard location (/usr/libexec/upowerd), a masquerading technique used by the Linux malware from the JDownloader compromise.'
references:
  - https://securityaffairs.co/wordpress/162985/malware/jdownloader-site-malware-supply-chain-attack.html
author: cyber.netsecops.io
date: 2026-05-10
tags:
  - attack.t1036.005
  - attack.defense_evasion
level: high
logsource:
  category: process_creation
  product: linux
detection:
  selection:
    Image: '/usr/libexec/upowerd'
  filter_legit:
    ParentImage: '/usr/lib/systemd/systemd'
  condition: selection and not filter_legit