# sigma/oxloader-uses-malicious-google-ads-to-spread-castlestealer-SIGMA.yaml
title: Suspicious PowerShell Execution Policy Bypass For File Execution
id: 6c8d7b6a-5e4f-4d3c-b2a1-9e8f7d6c5b4a
status: experimental
description: Detects PowerShell being launched with the ExecutionPolicy Bypass flag to run a script file, a technique often used by malware to run untrusted code.
references:
  - https://thehackernews.com/2026/06/new-oxloader-loader-uses-malicious.html
author: cyber.netsecops.io
date: 2026-06-22
tags:
  - attack.t1059.001
  - attack.execution
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    CommandLine|contains|all:
      - 'powershell'
      - '-ExecutionPolicy Bypass'
      - '-File'
  condition: selection
---
title: Suspicious Execution from User Temp Directory
id: 7d6c5b4a-8e7f-4c3b-a19e-8f7d6c5b4a9e
status: experimental
description: Detects processes being executed from common user temporary or application data directories, which is a common tactic for malware droppers like OXLOADER.
references:
  - https://thehackernews.com/2026/06/new-oxloader-loader-uses-malicious.html
author: cyber.netsecops.io
date: 2026-06-22
tags:
  - attack.t1204.002
  - attack.execution
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|contains:
      - ':\Users\Public\'
      - '\AppData\Local\Temp\'
      - '\AppData\Roaming\'
  filter_legit:
    Image|endswith:
      - '\Updater.exe'
      - '\setup.exe'
  condition: selection and not filter_legit
---
title: Suspicious Modification of Registry Run Key
id: 8e7f6d5c-9f8e-4b2a-9e8f-7d6c5b4a9e8f
status: experimental
description: Detects the creation or modification of values in the user's Run key, a common persistence technique used by malware.
references:
  - https://thehackernews.com/2026/06/new-oxloader-loader-uses-malicious.html
author: cyber.netsecops.io
date: 2026-06-22
tags:
  - attack.t1547.001
  - attack.persistence
level: medium
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    TargetObject|contains: '\Software\Microsoft\Windows\CurrentVersion\Run'
  condition: selection
---
title: OXLOADER Malicious and Staging Domain DNS Query
id: 5b9e1a2f-3c4d-4e5f-8a9b-0c1d2e3f4a5b
status: experimental
description: Detects DNS queries to domains associated with the OXLOADER malvertising campaign, including the fraudulent download site and the decentralized storage gateway.
references:
  - https://thehackernews.com/2026/06/new-oxloader-loader-uses-malicious.html
author: cyber.netsecops.io
date: 2026-06-22
tags:
  - attack.t1105
  - attack.command_and_control
level: high
logsource:
  category: dns_query
  product: windows
detection:
  selection:
    QueryName|contains:
      - 'node-js.prentiva99.info'
      - 'gateway.storjshare.io'
  condition: selection