# sigma/sidecopy-apt-targets-afghanistans-finance-ministry-with-xenorat-SIGMA.yaml
title: Suspicious Remote HTA Execution via Mshta.exe
id: 9a8b7c6d-5e4f-4321-a0b9-c8d7e6f5a4b3
status: experimental
description: Detects the use of mshta.exe to execute a remote HTML Application (HTA) file, a technique used by SideCopy in the XENOFISCAL campaign.
references:
  - https://thehackernews.com/2026/06/pakistan-linked-sidecopy-targets.html
author: cyber.netsecops.io
date: 2026-06-02
tags:
  - attack.t1218.005
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\mshta.exe'
    CommandLine|contains:
      - 'http://'
      - 'https://'
  condition: selection
---
title: Persistence via HKCU Run Key Creation
id: 1a2b3c4d-5e6f-4789-890a-bcdef0123456
status: experimental
description: 'Detects the creation of a new value in the HKCU\...\Run registry key, a common persistence technique used by malware like XenoRAT.'
references:
  - https://thehackernews.com/2026/06/pakistan-linked-sidecopy-targets.html
author: cyber.netsecops.io
date: 2026-06-02
tags:
  - attack.t1547.001
level: medium
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    EventType: 'SetValue'
    TargetObject|contains: '\Software\Microsoft\Windows\CurrentVersion\Run'
  condition: selection