# sigma/stealthy-powmix-botnet-targets-czech-republic-workforce-SIGMA.yaml
title: PowMix PowerShell Loader Execution
id: 9c5b1a3e-f8d2-4e1c-b9a5-6f7d8e9c0b1a
status: experimental
description: Detects the execution of a PowerShell script with a bypassed execution policy, a technique used by PowMix's LNK-based loader.
references:
  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html
author: cyber.netsecops.io
date: 2026-04-17
tags:
  - attack.t1059.001
  - attack.execution
level: high
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\powershell.exe'
    CommandLine|contains|all:
      - '-ExecutionPolicy Bypass'
      - '-File'
  condition: selection
---
title: PowMix Persistence via Scheduled Task
id: 3d2e1c0b-a9b8-4f7d-8e6c-5b4a3d2e1c0b
status: experimental
description: Detects the creation of a new scheduled task using schtasks.exe, which is the persistence mechanism used by the PowMix botnet.
references:
  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html
author: cyber.netsecops.io
date: 2026-04-17
tags:
  - attack.t1053.005
  - attack.persistence
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/Create'
  condition: selection
---
title: PowMix C2 URL Pattern
id: 7a6b5c4d-3e2f-4d1c-a9b8-c7d6e5f4a3b2
status: experimental
description: Detects outbound web requests with URL paths matching the pattern used by PowMix for C2 communications, which embeds data to mimic REST API calls.
references:
  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html
author: cyber.netsecops.io
date: 2026-04-17
tags:
  - attack.t1071.001
  - attack.command_and_control
level: medium
logsource:
  category: proxy
detection:
  selection:
    c-uri|re: '/api/v\d+/data=.*'
  condition: selection