# sigma/uac-0247-espionage-campaign-targets-government-and-healthcare-in-ukraine-SIGMA.yaml
title: Suspicious HTA File Execution by Email Client or Browser
id: f1b2c3d4-e5f6-4a7b-8c9d-0e1f2a3b4c5d
status: experimental
description: Detects mshta.exe being launched by a common email client or web browser, which is a TTP used by UAC-0247.
references:
  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html
author: cyber.netsecops.io
date: 2026-04-18
tags:
  - attack.t1218.005
  - attack.execution
level: high
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\mshta.exe'
    ParentImage|endswith:
      - '\outlook.exe'
      - '\chrome.exe'
      - '\firefox.exe'
      - '\msedge.exe'
  condition: selection
---
title: Suspicious Scheduled Task Creation by Script Interpreter
id: a2b3c4d5-e6f7-4b8c-9d0e-1f2a3b4c5d6e
status: experimental
description: Detects the creation of a scheduled task by a script interpreter like mshta.exe, a persistence technique used by UAC-0247.
references:
  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html
author: cyber.netsecops.io
date: 2026-04-18
tags:
  - attack.t1053.005
  - attack.persistence
level: medium
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
    ParentImage|endswith:
      - '\mshta.exe'
      - '\powershell.exe'
      - '\wscript.exe'
      - '\cscript.exe'
  condition: selection
---
title: Anomalous RuntimeBroker Network Activity
id: b3c4d5e6-f7a8-4c9d-0e1f-2a3b4c5d6f7a
status: experimental
description: Detects RuntimeBroker.exe initiating outbound network connections, which is atypical behavior and could indicate process injection as used by UAC-0247.
references:
  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html
author: cyber.netsecops.io
date: 2026-04-18
tags:
  - attack.t1055
  - attack.defense_evasion
level: medium
logsource:
  category: network_connection
  product: windows
detection:
  selection:
    Image|endswith: '\RuntimeBroker.exe'
    Initiated: 'true'
  filter:
    DestinationIp|cidr:
      - '127.0.0.0/8'
      - '::1/128'
  condition: selection and not filter