# data/articles/bloody-wolf-apt-targets-russia-and-uzbekistan-with-netsupport-rat.json
{
  "id": "03e306c5-5b5c-4185-8b95-ad8fe979a97f",
  "slug": "bloody-wolf-apt-targets-russia-and-uzbekistan-with-netsupport-rat",
  "headline": "'Bloody Wolf' APT Deploys NetSupport RAT in Espionage Campaign",
  "title": "'Bloody Wolf' APT Targets Russia and Uzbekistan with NetSupport RAT in Financial and Espionage Attacks",
  "summary": "Security researchers have uncovered an active spear-phishing campaign attributed to the threat actor 'Bloody Wolf' (also tracked as Stan Ghouls). The campaign targets organizations primarily in Uzbekistan and Russia, with a focus on manufacturing, finance, and IT sectors, though government and other entities have also been targeted. The attackers use phishing emails with password-protected ZIP archives containing a malicious LNK file. When executed, this file downloads and installs the legitimate remote administration tool, NetSupport RAT, which gives the attackers full control over the victim's system. The motives appear to be mixed, pointing towards both financially motivated cybercrime and state-aligned cyber espionage. This campaign marks a shift in tooling for the group, which previously used the STRRAT malware.",
  "full_report": "## Executive Summary\nAn ongoing spear-phishing campaign by the threat actor known as **Bloody Wolf** (tracked by Kaspersky as Stan Ghouls) is targeting organizations in Uzbekistan and Russia. The campaign, which has compromised approximately 60 victims, leverages the legitimate remote access tool **[NetSupport RAT](https://www.netsupportsoftware.com/)** to gain full control over infected systems. The attackers' targeting is broad, encompassing manufacturing, finance, IT, and government sectors, suggesting a dual motive of financial gain and cyber espionage. The attack begins with a phishing email containing a malicious LNK file disguised as a PDF within a password-protected ZIP archive. This campaign represents a tactical evolution for Bloody Wolf, which has shifted from using STRRAT malware to abusing a legitimate commercial tool for its operations, a technique known as Living off the Land (LotL).\n\n## Threat Overview\n- **Threat Actor:** Bloody Wolf (also known as Stan Ghouls).\n- **Activity Period:** Active since at least 2023.\n- **Targets:** Approximately 60 victims, with 50 in Uzbekistan and 10 in Russia. Targeted sectors include manufacturing, finance, IT, government, logistics, healthcare, and education.\n- **Motive:** Appears to be a hybrid of financial crime (targeting financial institutions) and cyber espionage (broad targeting and use of a RAT).\n- **Payload:** NetSupport RAT, a commercially available remote administration tool.\n- **Initial Access Vector:** Spear-phishing emails containing a password-protected ZIP archive. ([`T1566.001 - Spearphishing Attachment`](https://attack.mitre.org/techniques/T1566/001/))\n\n## Technical Analysis\nThe attack chain is straightforward and relies on tricking the user into executing the initial payload.\n1.  **Delivery:** The victim receives a phishing email with a password-protected ZIP file. Including the password in the email body is a common tactic to bypass automated email gateway scanners, which cannot inspect the archive's contents.\n2.  **Lure:** Inside the ZIP archive is a Windows shortcut (LNK) file that is disguised with a PDF icon and a deceptive name to trick the user into double-clicking it. ([`T1204.002 - Malicious File`](https://attack.mitre.org/techniques/T1204/002/))\n3.  **Execution:** When the LNK file is executed, it runs a command to download a loader from an external domain. ([`T1218.005 - Msiexec`](https://attack.mitre.org/techniques/T1218/005/) or similar LOLBins are often used by LNK files).\n4.  **Payload Installation:** The loader checks if NetSupport RAT is already installed. If not, it downloads the RAT from an attacker-controlled server and executes it.\n5.  **Persistence:** To ensure it survives a reboot, the malware creates an autorun script in the Windows Startup folder, which re-launches the RAT every time the user logs in. ([`T1547.001 - Registry Run Keys / Startup Folder`](https://attack.mitre.org/techniques/T1547/001/))\n6.  **Command and Control:** Once installed, the NetSupport RAT connects to the attacker's C2 server, giving them complete remote control over the victim's machine. They can then perform actions such as file transfer, keystroke logging, screen viewing, and executing further commands.\n\n> The use of a legitimate tool like NetSupport RAT is a classic Living off the Land, Binaries, and Scripts (LOLBAS) technique. It allows the attackers' activity to blend in with normal administrative traffic, making detection more difficult for security tools that rely on blacklisting known malicious files.\n\n## Impact Assessment\n- **Espionage:** With full remote access, the attackers can steal sensitive corporate data, intellectual property, government documents, and financial information.\n- **Financial Theft:** The focus on financial institutions suggests the attackers may use their access to commit fraud, such as initiating unauthorized wire transfers.\n- **Further Compromise:** The compromised machine can be used as a pivot point to move laterally within the victim's network, leading to a much wider breach.\n- **Installation of Other Malware:** The RAT provides a channel for the attackers to deploy additional malware, such as ransomware or banking trojans.\n\n## Detection & Response\n- **Network Traffic Monitoring:** Monitor for network traffic associated with NetSupport RAT, especially connections from general user workstations to external IP addresses on the default NetSupport ports (e.g., TCP 5405). Since it's a legitimate tool, context is key. An IT admin using it is normal; a user in the accounting department having it run is not.\n- **Endpoint Detection:** Use an EDR to alert on the installation and execution of `client32.exe`, the main executable for NetSupport. Look for the specific process chain of a LNK file leading to the download and execution of a new program.\n- **Log Analysis:** Analyze logs for the creation of files or scripts in the Startup folder (`%APPDATA%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup`).\n- **Email Security:** While scanning password-protected ZIPs is difficult, email gateways can be configured to flag or block such attachments for manual review, especially from untrusted senders.\n\n## Mitigation\n- **Application Control:** Implement application allow-listing to prevent unauthorized software like NetSupport RAT from being installed or executed on user endpoints. If NetSupport is not a sanctioned tool in your environment, it should be explicitly blocked.\n- **User Training:** Train users to be suspicious of emails with password-protected attachments and to never execute files from untrusted sources, especially LNK files disguised as documents.\n- **Attack Surface Reduction:** Configure Windows to show file extensions by default. This makes it easier for users to see that a file is a `.lnk` (a shortcut) and not a `.pdf` (a document).\n- **PowerShell Hardening:** Use constrained language mode and script block logging for PowerShell to limit the capabilities of malicious scripts and improve visibility into their execution.",
  "twitter_post": "The 'Bloody Wolf' APT is targeting orgs in Russia & Uzbekistan with a spear-phishing campaign that deploys the legitimate NetSupport RAT. The group's motives appear to be a mix of espionage and financial crime. ðŸº #APT #CyberAttack #NetSupport",
  "meta_description": "The 'Bloody Wolf' (Stan Ghouls) threat actor is conducting a spear-phishing campaign against entities in Russia and Uzbekistan, using the NetSupport RAT for espionage and financial gain.",
  "category": [
    "Threat Actor",
    "Cyberattack",
    "Malware"
  ],
  "severity": "high",
  "entities": [
    {
      "name": "Bloody Wolf",
      "type": "threat_actor"
    },
    {
      "name": "NetSupport RAT",
      "type": "malware"
    },
    {
      "name": "Kaspersky",
      "type": "vendor"
    },
    {
      "name": "STRRAT",
      "type": "malware"
    }
  ],
  "cves": [],
  "sources": [
    {
      "url": "https://thehackernews.com/2026/02/bloody-wolf-targets-uzbekistan-russia.html",
      "title": "Bloody Wolf Targets Uzbekistan, Russia Using NetSupport RAT in Spear-Phishing Campaign",
      "date": "2026-02-09",
      "friendly_name": "The Hacker News",
      "website": "thehackernews.com"
    },
    {
      "url": "https://www.wiu.edu/cbt/cybersecurity/news.php",
      "title": "Cybersecurity News",
      "date": "2026-02-09",
      "friendly_name": "WIU Cybersecurity",
      "website": "wiu.edu"
    }
  ],
  "events": [],
  "mitre_techniques": [
    {
      "id": "T1566.001",
      "name": "Spearphishing Attachment",
      "tactic": "Initial Access"
    },
    {
      "id": "T1204.002",
      "name": "Malicious File",
      "tactic": "Execution"
    },
    {
      "id": "T1547.001",
      "name": "Registry Run Keys / Startup Folder",
      "tactic": "Persistence"
    },
    {
      "id": "T1219",
      "name": "Remote Access Software",
      "tactic": "Command and Control"
    },
    {
      "id": "T1053.005",
      "name": "Scheduled Task",
      "tactic": "Persistence"
    }
  ],
  "mitre_mitigations": [
    {
      "id": "M1038",
      "name": "Execution Prevention",
      "d3fend_techniques": [
        {
          "id": "D3-EAL",
          "name": "Executable Allowlisting",
          "url": "https://d3fend.mitre.org/technique/d3f:ExecutableAllowlisting"
        }
      ],
      "description": "Application allow-listing can prevent the execution of unauthorized remote access tools like NetSupport RAT if it is not approved for use in the environment.",
      "domain": "enterprise"
    },
    {
      "id": "M1017",
      "name": "User Training",
      "description": "Training users to recognize phishing attempts, especially those with suspicious attachments like password-protected ZIPs and LNK files, is a key preventative measure.",
      "domain": "enterprise"
    },
    {
      "id": "M1037",
      "name": "Filter Network Traffic",
      "d3fend_techniques": [
        {
          "id": "D3-OTF",
          "name": "Outbound Traffic Filtering",
          "url": "https://d3fend.mitre.org/technique/d3f:OutboundTrafficFiltering"
        }
      ],
      "description": "Blocking outbound connections on ports used by remote access tools like NetSupport RAT (e.g., TCP 5405) can prevent the malware from connecting to its C2 server.",
      "domain": "enterprise"
    }
  ],
  "d3fend_countermeasures": [
    {
      "technique_id": "D3-EAL",
      "technique_name": "Executable Allowlisting",
      "url": "https://d3fend.mitre.org/technique/d3f:ExecutableAllowlisting",
      "recommendation": "To counter the 'Bloody Wolf' campaign's abuse of a legitimate tool, executable allowlisting is a highly effective control. Since NetSupport RAT is commercial software, its executables (like `client32.exe`) are signed and may not be flagged as malicious by traditional antivirus. An allowlisting approach, however, would block its execution by default on all systems where it is not explicitly authorized. Security teams should create a policy that only permits known, approved applications to run on user endpoints. If NetSupport is not used for legitimate IT administration, it should not be on the allowlist. This prevents the payload from running, even if a user is tricked into executing the initial LNK file. This D3FEND technique effectively neutralizes the attacker's 'Living off the Land' strategy by treating any unauthorized legitimate tool as malicious.",
      "mitre_mitigation_id": "M1038"
    },
    {
      "technique_id": "D3-PA",
      "technique_name": "Process Analysis",
      "url": "https://d3fend.mitre.org/technique/d3f:ProcessAnalysis",
      "recommendation": "Detecting this attack requires analyzing process lineage. An EDR solution should be configured to alert on suspicious parent-child process relationships. In this campaign, the chain would look something like `explorer.exe` (user double-clicks LNK) -> `cmd.exe` or `powershell.exe` (executed by the LNK file) -> network connection to download the loader -> execution of the NetSupport RAT installer. This is an anomalous chain of events. A user's file explorer should not be leading to the installation of a remote access tool. By creating detection rules that look for this specific sequence of process creation and network activity, security teams can identify the attack at the execution stage, before the RAT is fully installed and connected to its C2 server. This provides a critical window for automated or manual intervention.",
      "mitre_mitigation_id": "M1049"
    }
  ],
  "iocs": [],
  "cyber_observables": [
    {
      "type": "file_name",
      "value": "*.lnk",
      "description": "The initial payload is a Windows shortcut (LNK) file. The presence of LNK files within ZIP archives from external emails is highly suspicious.",
      "context": "Email gateway attachment analysis, Endpoint file monitoring",
      "confidence": "medium"
    },
    {
      "type": "process_name",
      "value": "client32.exe",
      "description": "The main process name for the NetSupport Manager client. Its execution on non-IT administrator workstations is a strong indicator of compromise.",
      "context": "EDR, Process monitoring",
      "confidence": "high"
    },
    {
      "type": "file_path",
      "value": "%APPDATA%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup",
      "description": "The location used by the malware to place a script for persistence. Monitor for new file creation in this directory.",
      "context": "File integrity monitoring, EDR",
      "confidence": "high"
    },
    {
      "type": "port",
      "value": "5405",
      "description": "The default TCP port used by NetSupport Manager for communication. Outbound connections on this port from user endpoints should be investigated.",
      "context": "Firewall logs, NetFlow analysis",
      "confidence": "medium"
    }
  ],
  "tags": [
    "Bloody Wolf",
    "APT",
    "NetSupport RAT",
    "Spear-phishing",
    "Uzbekistan",
    "Russia",
    "LOLBAS"
  ],
  "extract_datetime": "2026-02-09T15:00:00.000Z",
  "article_type": "NewsArticle",
  "impact_scope": {
    "geographic_scope": "regional",
    "countries_affected": [
      "Uzbekistan",
      "Russia",
      "Kyrgyzstan",
      "Kazakhstan"
    ],
    "industries_affected": [
      "Manufacturing",
      "Finance",
      "Technology",
      "Government",
      "Healthcare",
      "Education"
    ]
  },
  "pub_date": "2026-02-09",
  "reading_time_minutes": 5,
  "createdAt": "2026-02-09T15:00:00.000Z",
  "updatedAt": "2026-02-09T15:00:00.000Z"
}