# data/articles/google-disrupts-sprawling-chinese-espionage-campaign-unc2814.json
{
  "id": "c35c4bca-3bf6-4b50-83d9-149d3eced744",
  "slug": "google-disrupts-sprawling-chinese-espionage-campaign-unc2814",
  "headline": "Google & Mandiant Dismantle Global Chinese Spy Network Using Novel \"GRIDTIDE\" Backdoor",
  "title": "Google and Mandiant Disrupt Long-Running Chinese Espionage Campaign by UNC2814 Leveraging Google Sheets for C2",
  "summary": "On February 25, 2026, Google and Mandiant revealed they had disrupted a massive, multi-year cyber espionage campaign attributed to UNC2814, a suspected China-nexus threat actor. The operation compromised at least 53 organizations in 42 countries, primarily in the telecommunications and government sectors. The attackers utilized a novel C-based backdoor named GRIDTIDE, which cleverly abused the legitimate Google Sheets API for command and control, allowing its traffic to blend in with normal cloud activity. The joint disruption effort involved terminating the actor's cloud projects, disabling accounts, and taking down associated infrastructure, marking a significant blow to a prolific and distinct espionage group.",
  "full_report": "## Executive Summary\nOn February 25, 2026, **[Google](https://cloud.google.com/blog/threat-intelligence/)** and **[Mandiant](https://www.mandiant.com/)** announced the successful disruption of a global cyber espionage campaign conducted by **[UNC2814](https://www.mandiant.com/resources/blog/unc2814-espionage-campaign)**, a sophisticated threat actor with suspected links to the People's Republic of China (PRC). The campaign, active since at least 2017, utilized a novel backdoor named **GRIDTIDE** that employed the **[Google Sheets API](https://developers.google.com/sheets/api)** for command and control (C2), enabling it to evade detection for years. The operation compromised at least 53 organizations across 42 countries, with a primary focus on telecommunications and government entities. The joint takedown involved a coordinated effort to dismantle the actor's infrastructure, though UNC2814 is expected to attempt to rebuild its capabilities.\n\n## Threat Overview\nThe threat actor, tracked as **UNC2814**, has demonstrated a high level of operational security and sophistication, with no observed overlap with other known Chinese APT groups. The campaign's primary goal was long-term intelligence gathering. The main targets were government and telecommunications organizations in Africa, Asia, and the Americas, suggesting a strategic focus on geopolitical and economic intelligence.\n\nThe core of the operation was the **GRIDTIDE** backdoor, a custom malware written in C. Its most innovative feature was its C2 mechanism. By using the legitimate Google Sheets API, attackers could issue commands and exfiltrate data without generating suspicious network traffic that would typically be flagged by security solutions. This living-off-the-trusted-land technique highlights a growing trend of adversaries abusing legitimate cloud services to mask their activities.\n\n## Technical Analysis\n**Initial Access:** While the specific initial access vector for all intrusions remains under investigation, **UNC2814** has a known history of exploiting vulnerabilities in public-facing web servers and other edge network systems, a common tactic mapped to [`T1190 - Exploit Public-Facing Application`](https://attack.mitre.org/techniques/T1190/).\n\n**Execution & Persistence:** Once inside a network, the attackers established persistence for the **GRIDTIDE** backdoor by creating a systemd service. This was achieved by creating a service file at `/etc/systemd/system/xapt.service`, a technique identified as [`T1543.002 - Create or Modify System Process: Systemd Service`](https://attack.mitre.org/techniques/T1543/002/). The backdoor itself would then be executed with system-level privileges.\n\n**Command and Control:** The C2 mechanism is the most notable aspect of this campaign. **GRIDTIDE** communicated via the Google Sheets API, a form of [`T1102 - Web Service`](https://attack.mitre.org/techniques/T1102/).\n1.  The malware would periodically poll cell `A1` of a specific Google Sheet for new commands.\n2.  Victim metadata, such as hostname and user information, was written to cell `V1`.\n3.  Other cells were used for staging data for exfiltration and receiving command output.\n4.  All data transferred via the API was encoded using a URL-safe Base64 scheme to avoid detection by web filters.\n\n**Lateral Movement:** The threat actor used service accounts to move laterally within compromised networks via [`T1021.004 - Remote Services: SSH`](https://attack.mitre.org/techniques/T1021/004/). This allowed them to pivot from the initial point of entry to other high-value systems while using legitimate credentials and protocols to evade detection.\n\n## Impact Assessment\nThe impact of this campaign is significant due to its longevity, scale, and the nature of the targeted organizations. By compromising telecommunications providers and government agencies for years, **UNC2814** likely gained access to vast amounts of sensitive data, including communications records, government documents, and intellectual property. The long-term persistence within these networks suggests a deeply entrenched intelligence-gathering operation. For the 53 confirmed victim organizations, the breach represents a severe compromise of security and trust, necessitating extensive and costly incident response, forensic analysis, and system remediation efforts to ensure the complete eviction of the threat actor.\n\n## Cyber Observables for Detection\nSecurity teams should hunt for the following patterns to detect potential **GRIDTIDE** or similar activity:\n\n| Type | Value | Description |\n|---|---|---|\n| API Endpoint | `sheets.googleapis.com` | Monitor for anomalous or high-frequency API calls from servers, especially those not expected to interact with Google Sheets. |\n| User Agent | `Google-API-C++` | While potentially legitimate, look for this user agent from unexpected sources or in conjunction with other suspicious indicators. |\n| File Path | `/etc/systemd/system/xapt.service` | Monitor for the creation of this specific service file or other new, unrecognized systemd services. |\n| Network Traffic Pattern | Outbound connections to `sheets.googleapis.com` on port 443 | Scrutinize TLS traffic to this endpoint from servers, focusing on unusual data volumes or connection frequencies. |\n\n## Detection & Response\nDetecting this type of activity requires a multi-layered approach that goes beyond traditional network signatures.\n\n1.  **Cloud API Monitoring:** Implement robust monitoring of all API traffic to and from cloud services like Google Workspace. Establish baselines for normal API usage and alert on deviations, such as servers unexpectedly communicating with the Google Sheets API. This aligns with the D3FEND technique `Cloud API Monitoring`.\n2.  **Endpoint Detection and Response (EDR):** Deploy EDR solutions to monitor for suspicious process creation and persistence mechanisms. A rule to detect the creation of new systemd services (`/etc/systemd/system/`) can be highly effective. This relates to D3FEND's [`D3-PA - Process Analysis`](https://d3fend.mitre.org/technique/d3f:ProcessAnalysis).\n3.  **Log Analysis:** Correlate network logs with endpoint logs. Look for SSH logins using service accounts originating from unusual internal sources. Analyze web proxy and DNS logs for connections to `sheets.googleapis.com` from non-standard clients.\n\n## Mitigation\nOrganizations should implement the following controls to defend against this and similar threats:\n\n1.  **Restrict Outbound Traffic:** Implement egress filtering to block or restrict access to cloud services and APIs that are not required for business operations. Deny access to personal Google accounts and other non-corporate cloud services from server segments. This is a form of D3FEND's [`D3-OTF - Outbound Traffic Filtering`](https://d3fend.mitre.org/technique/d3f:OutboundTrafficFiltering).\n2.  **Privileged Access Management (PAM):** Strictly control the use of service accounts. Ensure they have the minimum necessary privileges and are monitored for anomalous activity. Limit where service accounts can log in from and disable interactive logon capabilities where possible.\n3.  **Application Control:** Use application control solutions to prevent unauthorized software, including custom backdoors, from executing on servers.\n4.  **Patch Management:** Maintain a rigorous patch management program to remediate vulnerabilities in public-facing applications, which are common initial access vectors for actors like **UNC2814**.",
  "twitter_post": "Google & Mandiant disrupt massive Chinese espionage op by UNC2814! ⚠️ Attackers hit 53 orgs in 42 countries using novel \"GRIDTIDE\" backdoor that abuses Google Sheets for C2. #CyberEspionage #ThreatIntel #UNC2814",
  "meta_description": "Google and Mandiant have taken down a global cyber espionage campaign by UNC2814, a China-linked APT that used the GRIDTIDE backdoor and Google Sheets for C2 to breach 53 organizations.",
  "category": [
    "Threat Actor",
    "Cyberattack",
    "Threat Intelligence"
  ],
  "severity": "high",
  "entities": [
    {
      "name": "UNC2814",
      "type": "threat_actor"
    },
    {
      "name": "GRIDTIDE",
      "type": "malware"
    },
    {
      "name": "Google",
      "type": "vendor",
      "url": "https://cloud.google.com/blog/threat-intelligence/"
    },
    {
      "name": "Mandiant",
      "type": "security_organization",
      "url": "https://www.mandiant.com/"
    },
    {
      "name": "Google Sheets API",
      "type": "technology"
    }
  ],
  "cves": [],
  "sources": [
    {
      "url": "https://cloud.google.com/blog/threat-intelligence/disrupting-unc2814-global-espionage-campaign",
      "title": "Disrupting the GRIDTIDE Global Cyber Espionage Campaign",
      "date": "2026-02-25",
      "friendly_name": "Google Cloud Blog",
      "website": "cloud.google.com"
    },
    {
      "url": "https://thehackernews.com/2026/02/google-disrupts-unc2814-gridtide.html",
      "title": "Google Disrupts UNC2814 GRIDTIDE Campaign After 53 Breaches Across 42 Countries",
      "date": "2026-02-25",
      "friendly_name": "The Hacker News",
      "website": "thehackernews.com"
    }
  ],
  "events": [
    {
      "datetime": "2017-01-01T00:00:00Z",
      "summary": "The UNC2814 espionage campaign is believed to have started its operations."
    },
    {
      "datetime": "2026-02-25T00:00:00Z",
      "summary": "Google and Mandiant publicly announce the disruption of the UNC2814 campaign and its GRIDTIDE backdoor."
    }
  ],
  "mitre_techniques": [
    {
      "id": "T1102",
      "name": "Web Service",
      "tactic": "Command and Control"
    },
    {
      "id": "T1543.002",
      "name": "Create or Modify System Process: Systemd Service",
      "tactic": "Persistence"
    },
    {
      "id": "T1021.004",
      "name": "Remote Services: SSH",
      "tactic": "Lateral Movement"
    },
    {
      "id": "T1190",
      "name": "Exploit Public-Facing Application",
      "tactic": "Initial Access"
    },
    {
      "id": "T1059.004",
      "name": "Unix Shell",
      "tactic": "Execution"
    }
  ],
  "mitre_mitigations": [
    {
      "id": "M1037",
      "name": "Filter Network Traffic",
      "domain": "enterprise"
    },
    {
      "id": "M1021",
      "name": "Restrict Web-Based Content",
      "description": "Restrict access to cloud APIs and services from server environments where it is not explicitly required for business operations.",
      "domain": "enterprise"
    },
    {
      "id": "M1026",
      "name": "Privileged Account Management",
      "description": "Monitor and control the usage of service accounts, especially for interactive logons and remote access like SSH.",
      "domain": "enterprise"
    },
    {
      "id": "M1047",
      "name": "Audit",
      "description": "Enable and review logs for cloud API usage, endpoint process creation, and remote service authentications.",
      "domain": "enterprise"
    }
  ],
  "d3fend_countermeasures": [
    {
      "technique_id": "D3-OTF",
      "technique_name": "Outbound Traffic Filtering",
      "url": "https://d3fend.mitre.org/technique/d3f:OutboundTrafficFiltering",
      "recommendation": "To counter the GRIDTIDE backdoor's C2 method, implement strict outbound traffic filtering on server subnets. Specifically, create firewall rules that deny all outbound connections to `sheets.googleapis.com` and other non-essential Google API endpoints by default. An allowlist should be created for specific servers that have a legitimate business need to access these services. This forces adversary traffic into monitored and controlled channels, preventing backdoors like GRIDTIDE from establishing C2. This should be coupled with TLS inspection where possible to analyze the content of the API calls, further enhancing detection capabilities for malicious use of legitimate services.",
      "mitre_mitigation_id": "M1037"
    },
    {
      "technique_id": "D3-PA",
      "technique_name": "Process Analysis",
      "url": "https://d3fend.mitre.org/technique/d3f:ProcessAnalysis",
      "recommendation": "To detect GRIDTIDE's persistence mechanism, deploy an EDR solution capable of deep process analysis on Linux servers. Configure detection rules to specifically monitor for and alert on the creation of new files within the `/etc/systemd/system/` directory. Establish a baseline of known-good services and alert on any deviation. Furthermore, monitor for processes that make network connections to `sheets.googleapis.com`. A server process that has no legitimate reason to interact with Google Sheets suddenly initiating such connections is a strong indicator of compromise. This proactive monitoring directly targets the TTPs used by UNC2814 to establish a foothold.",
      "mitre_mitigation_id": "M1049"
    },
    {
      "technique_id": "D3-LAM",
      "technique_name": "Local Account Monitoring",
      "url": "https://d3fend.mitre.org/technique/d3f:LocalAccountMonitoring",
      "recommendation": "To address UNC2814's lateral movement via SSH with service accounts, enhance account monitoring. Implement SIEM rules to alert on any interactive logon using a service account. These accounts should typically only be used by applications and services, not for remote shell access. Correlate SSH authentication logs (e.g., from `/var/log/auth.log` on Linux) with a list of designated service accounts. Any successful SSH login by an account on this list should trigger a high-priority security alert. This helps detect the misuse of legitimate credentials for lateral movement, a key part of the UNC2814 playbook.",
      "mitre_mitigation_id": "M1026"
    }
  ],
  "iocs": [],
  "cyber_observables": [
    {
      "type": "api_endpoint",
      "value": "sheets.googleapis.com/v4/spreadsheets/",
      "description": "The API endpoint used by GRIDTIDE for command and control via Google Sheets.",
      "context": "Network traffic logs, cloud audit logs. Monitor for anomalous requests from servers.",
      "confidence": "high"
    },
    {
      "type": "file_path",
      "value": "/etc/systemd/system/xapt.service",
      "description": "Persistence mechanism used by the GRIDTIDE backdoor to run as a service.",
      "context": "File integrity monitoring (FIM), EDR logs, or manual inspection on Linux systems.",
      "confidence": "high"
    },
    {
      "type": "log_source",
      "value": "Google Workspace Audit Logs",
      "description": "Logs showing API calls to Google Sheets can reveal unauthorized access patterns.",
      "context": "SIEM or cloud security monitoring platform. Look for activity from service accounts or server IP addresses.",
      "confidence": "medium"
    },
    {
      "type": "command_line_pattern",
      "value": "ssh -l <service_account>",
      "description": "UNC2814 used service accounts for lateral movement via SSH.",
      "context": "Monitor command line logs (Event ID 4688, auditd) for SSH connections using service account credentials.",
      "confidence": "medium"
    }
  ],
  "tags": [
    "UNC2814",
    "GRIDTIDE",
    "Cyber Espionage",
    "APT",
    "Google Sheets",
    "C2",
    "Threat Intelligence",
    "China"
  ],
  "extract_datetime": "2026-02-25T15:00:00.000Z",
  "article_type": "NewsArticle",
  "impact_scope": {
    "geographic_scope": "global",
    "countries_affected": [
      "42 countries across Africa, Asia, and the Americas"
    ],
    "industries_affected": [
      "Telecommunications",
      "Government"
    ],
    "other_affected": [
      "Suspected infections in an additional 20 nations"
    ]
  },
  "pub_date": "2026-02-25",
  "reading_time_minutes": 5,
  "createdAt": "2026-02-25T15:00:00.000Z",
  "updatedAt": "2026-02-25T15:00:00.000Z"
}