# data/articles/huntress-uncovers-adware-campaign-exposing-25000-systems-to-supply-chain-attack.json
{
  "id": "6fcfbd76-a80a-4987-9e57-5c19a01e5f95",
  "slug": "huntress-uncovers-adware-campaign-exposing-25000-systems-to-supply-chain-attack",
  "headline": "Adware with Fangs: 25,000 Systems Exposed to $10 Supply Chain Hijack by Dragon Boss Solutions",
  "title": "Huntress Uncovers Adware from Dragon Boss Solutions That Disabled AV and Exposed 25,000+ Systems to Trivial Supply Chain Attack",
  "summary": "Security firm Huntress has exposed a dangerous operation where adware signed by 'Dragon Boss Solutions' went far beyond typical potentially unwanted programs (PUPs). The software, found on over 25,000 endpoints, used SYSTEM privileges to disable antivirus products and establish persistence. Critically, it used an unregistered domain for updates, `chromsterabrowser[.]com`, which could have been purchased for $10 by any attacker to push ransomware or other malware to all infected systems, including those in government, healthcare, and critical infrastructure networks. Huntress defensively registered the domain to prevent a widespread supply chain attack.",
  "full_report": "## Executive Summary\nResearchers at **[Huntress](https://www.huntress.com/blog)** have uncovered a high-risk threat masquerading as a common adware or Potentially Unwanted Program (PUP). The software, digitally signed by a UAE-based entity named **Dragon Boss Solutions**, was installed on at least 25,000 systems across 124 countries. While initially appearing as a simple browser hijacker, its capabilities evolved to include aggressive defense evasion techniques. Using a PowerShell payload with SYSTEM privileges, the software systematically disabled security products, blocked their update servers, and established persistence. The most alarming discovery was a critical flaw in its update mechanism: the software attempted to fetch updates from `chromsterabrowser[.]com`, a domain that was unregistered. This created a trivial but catastrophic supply chain risk, as any malicious actor could have registered the domain for about $10 and delivered malware to the entire botnet, which included universities, government agencies, and OT networks.\n\n---\n\n## Threat Overview\nThe operation demonstrates a dangerous evolution from adware to a potent backdoor. The software, once installed, executed a series of malicious actions to entrench itself on the host and eliminate security controls.\n\n**Key TTPs:**\n1.  **Defense Evasion:** A PowerShell script with `SYSTEM` privileges was used to disable a wide range of cybersecurity products. It would terminate their processes, modify the `hosts` file to block communication with update and telemetry servers, and delete registry keys to prevent reinstallation.\n2.  **Persistence:** The malware established persistence through multiple methods, including Windows Management Instrumentation (WMI) event subscriptions and scheduled tasks, ensuring it would survive reboots and removal attempts.\n3.  **Supply Chain Vulnerability:** The core of the threat lay in its insecure update process. The hardcoded update domain, `chromsterabrowser[.]com`, was not registered by the developers. This is a classic example of a dangling domain/subdomain takeover vulnerability. An attacker could simply purchase the domain and configure it to serve a malicious payload in response to the update check-in requests from all 25,000+ infected hosts.\n\n## Impact Assessment\nHuntress's quick action in registering the domain and sinkholing the traffic prevented a potentially devastating attack. The scale of the infection was vast, with 23,565 unique IPs connecting to the sinkhole in just 24 hours. The compromised hosts were not just consumer devices; the analysis identified numerous high-value targets:\n- **221** Universities and colleges\n- **41** Operational Technology (OT) networks (including electric utilities)\n- **35** Government entities\n- **3** Healthcare organizations\n\nA successful supply chain hijack could have led to widespread ransomware deployment, data theft, or espionage across sensitive sectors. The fact that the adware had already disabled local security tools means that any follow-on attack would have had a very high chance of success.\n\n## IOCs\n| Type | Value | Description |\n|---|---|---|\n| domain | `chromsterabrowser[.]com` | Unregistered update domain used by the malware. Now sinkholed by Huntress. |\n| domain | `worldwidewebframework3[.]com` | Another C2 domain associated with the campaign. |\n\n## Detection and Response\n- **Check for IOCs:** Scan network logs for any connections to the domains listed above. Search file systems and registry for artifacts related to Dragon Boss Solutions software.\n- **Review Disabled Services:** On endpoints, check for disabled or non-functioning antivirus and EDR services. Investigate any unauthorized modifications to the `hosts` file (located at `C:\\Windows\\System32\\drivers\\etc\\hosts`).\n- **Hunt for Persistence:** Use tools to inspect WMI event subscriptions and scheduled tasks for suspicious entries created by the adware.\n- **Remove the PUP:** If the software is detected, a thorough removal is required, which may involve manual deletion of files, registry keys, and persistence mechanisms, followed by the re-installation of security tools.\n\n## Mitigation\n- **Application Allowlisting:** Implement application control policies to prevent the execution of unauthorized or untrusted software, including PUPs and adware.\n- **PowerShell Logging:** Enable enhanced PowerShell logging (Module Logging, Script Block Logging) to capture and analyze the execution of potentially malicious scripts.\n- **DNS Sinkholing:** Organizations can proactively sinkhole suspicious or known-bad domains at their own DNS resolvers to prevent connections.\n- **Supply Chain Scrutiny:** This incident serves as a reminder that even seemingly low-risk software can introduce significant supply chain vulnerabilities. Vet all software, even if it is digitally signed.",
  "twitter_post": "🤯 A $10 domain could have hijacked 25,000+ systems! Huntress found adware from 'Dragon Boss Solutions' disabling AV and using an unregistered update domain. OT, gov, & healthcare networks were at risk. #SupplyChain #CyberSecurity #Adware",
  "meta_description": "Huntress researchers detail a sophisticated adware campaign by Dragon Boss Solutions that exposed over 25,000 systems, including OT and government networks, to a simple supply chain attack via an unregistered update domain.",
  "category": [
    "Supply Chain Attack",
    "Malware",
    "Threat Intelligence"
  ],
  "severity": "critical",
  "entities": [
    {
      "name": "Huntress",
      "type": "security_organization",
      "url": "https://www.huntress.com/"
    },
    {
      "name": "Dragon Boss Solutions",
      "type": "threat_actor"
    },
    {
      "name": "PowerShell",
      "type": "technology",
      "url": "https://attack.mitre.org/techniques/T1059/001/"
    }
  ],
  "cves": [],
  "sources": [
    {
      "url": "https://www.securityweek.com/10-domain-could-have-handed-hackers-25k-endpoints-including-in-ot-and-gov-networks/",
      "title": "$10 Domain Could Have Handed Hackers 25k Endpoints, Including in OT and Gov Networks",
      "friendly_name": "SecurityWeek",
      "website": "securityweek.com"
    },
    {
      "url": "https://www.huntress.com/blog/pups-grow-fangs-dragon-boss-solutions-supply-chain-risk",
      "title": "When PUPs Grow Fangs: Dragon Boss Solutions' $10 Supply Chain Risk",
      "friendly_name": "Huntress",
      "website": "huntress.com"
    },
    {
      "url": "https://cybernews.com/security/trusted-adware-app-left-25000-systems-open-to-a-10-supply-chain-hijack/",
      "title": "Trusted adware app left 25,000+ systems open to a $10 supply-chain hijack",
      "friendly_name": "Cybernews",
      "website": "cybernews.com"
    },
    {
      "url": "https://www.infosecurity-magazine.com/news/adware-exposes-25000-endpoints/",
      "title": "Adware Campaign Exposes 25,000 Endpoints to Supply Chain Hijack",
      "friendly_name": "Infosecurity Magazine",
      "website": "infosecurity-magazine.com"
    }
  ],
  "events": [
    {
      "datetime": "2025-03-01T00:00:00Z",
      "summary": "Adware begins deploying PowerShell-based payload to disable security products."
    },
    {
      "datetime": "2026-04-14T00:00:00Z",
      "summary": "Huntress publishes their research on the Dragon Boss Solutions adware campaign."
    }
  ],
  "mitre_techniques": [
    {
      "id": "T1587.003",
      "name": "Digital Certificates",
      "tactic": "Resource Development"
    },
    {
      "id": "T1559.002",
      "name": "Dynamic-Link Library Hijacking",
      "tactic": "Defense Evasion"
    },
    {
      "id": "T1059.001",
      "name": "PowerShell",
      "tactic": "Execution"
    },
    {
      "id": "T1543.003",
      "name": "Windows Service",
      "tactic": "Persistence"
    },
    {
      "id": "T1562.001",
      "name": "Disable or Modify Tools",
      "tactic": "Defense Evasion"
    },
    {
      "id": "T1112",
      "name": "Modify Registry",
      "tactic": "Defense Evasion"
    }
  ],
  "mitre_mitigations": [
    {
      "id": "M1038",
      "name": "Execution Prevention",
      "description": "Use application allowlisting to prevent unauthorized software like this adware from running in the first place.",
      "domain": "enterprise"
    },
    {
      "id": "M1021",
      "name": "Restrict Web-Based Content",
      "description": "Use DNS filtering services to block connections to known malicious or untrusted domains, including those used by adware.",
      "domain": "enterprise"
    },
    {
      "id": "M1047",
      "name": "Audit",
      "description": "Enable and monitor PowerShell script block logging and process creation events to detect malicious script execution and defense evasion techniques.",
      "domain": "enterprise"
    }
  ],
  "d3fend_countermeasures": [
    {
      "technique_id": "D3-EAL",
      "technique_name": "Executable Allowlisting",
      "url": "https://d3fend.mitre.org/technique/d3f:ExecutableAllowlisting",
      "recommendation": "The most effective preventative measure against threats like the Dragon Boss Solutions adware is application allowlisting. By configuring systems to only run explicitly approved applications, organizations can block the initial execution of the PUP. This approach is particularly effective in sensitive environments like OT networks and government systems, which were targeted in this campaign. Instead of relying on blocklists that are always a step behind attackers, allowlisting creates a 'default deny' posture. Implementation should start with critical servers and fixed-function workstations where the software environment is stable. Tools like Windows AppLocker or third-party solutions can be used to enforce these policies based on publisher, path, or file hash.",
      "mitre_mitigation_id": "M1038"
    },
    {
      "technique_id": "D3-DNSDL",
      "technique_name": "DNS Denylisting",
      "url": "https://d3fend.mitre.org/technique/d3f:DNSDenylisting",
      "recommendation": "To counter the supply chain risk from the unregistered domain `chromsterabrowser[.]com`, organizations should implement robust DNS filtering. This involves using a DNS security service that blocks access to known malicious domains, phishing sites, and domains associated with PUPs/adware. Even if the adware gets onto a system, DNS filtering can sever its connection to command-and-control and update servers, rendering it inert. Organizations should configure their DNS resolvers to use a security-focused service and ensure that endpoints cannot bypass this control. This technique would have prevented the infected clients from ever reaching the malicious update server, had an attacker registered the domain.",
      "mitre_mitigation_id": "M1021"
    }
  ],
  "iocs": [
    {
      "type": "domain",
      "value": "chromsterabrowser[.]com",
      "description": "Unregistered update domain used by the adware, now sinkholed by Huntress."
    },
    {
      "type": "domain",
      "value": "worldwidewebframework3[.]com",
      "description": "Associated Command and Control domain."
    }
  ],
  "cyber_observables": [
    {
      "type": "file_path",
      "value": "C:\\Windows\\System32\\drivers\\etc\\hosts",
      "description": "The malware modifies this file to block security vendor domains. Monitor for unauthorized changes.",
      "context": "File Integrity Monitoring (FIM)",
      "confidence": "high"
    },
    {
      "type": "command_line_pattern",
      "value": "powershell.exe -ExecutionPolicy Bypass -Command *Set-MpPreference -DisableRealtimeMonitoring $true*",
      "description": "A command pattern used to disable Microsoft Defender's real-time monitoring via PowerShell.",
      "context": "Windows Event ID 4688 (Process Creation) with command line logging enabled",
      "confidence": "high"
    },
    {
      "type": "registry_key",
      "value": "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\",
      "description": "The malware creates scheduled tasks for persistence. Look for suspicious tasks not associated with legitimate software.",
      "context": "Registry monitoring tools or manual inspection",
      "confidence": "medium"
    },
    {
      "type": "log_source",
      "value": "WMI Activity (Event ID 5861)",
      "description": "The malware uses WMI event subscriptions for persistence. Monitoring for new WMI permanent event consumer creation can detect this activity.",
      "context": "Windows Event Log (Microsoft-Windows-WMI-Activity/Operational)",
      "confidence": "medium"
    }
  ],
  "tags": [
    "Adware",
    "Supply Chain Attack",
    "Huntress",
    "Dragon Boss Solutions",
    "PUP",
    "PowerShell",
    "Defense Evasion"
  ],
  "extract_datetime": "2026-04-15T15:00:00.000Z",
  "article_type": "NewsArticle",
  "impact_scope": {
    "geographic_scope": "global",
    "countries_affected": [
      "United States"
    ],
    "industries_affected": [
      "Education",
      "Government",
      "Healthcare",
      "Critical Infrastructure"
    ],
    "other_affected": [
      "Operational Technology (OT) networks"
    ],
    "people_affected_estimate": "25,000+ systems"
  },
  "pub_date": "2026-04-15",
  "reading_time_minutes": 4,
  "createdAt": "2026-04-15T15:00:00.000Z",
  "updatedAt": "2026-04-15T15:00:00.000Z"
}