# data/articles/new-baradai-ransomware-targets-multiple-industries-worldwide.json
{
  "id": "a5064bee-5acb-4d43-913f-5d1904f0c79c",
  "slug": "new-baradai-ransomware-targets-multiple-industries-worldwide",
  "headline": "New 'BARADAI' Ransomware Emerges, Demanding Up to $80,000",
  "title": "CYFIRMA Discovers New BARADAI Ransomware Strain Targeting Windows Systems Across Multiple Sectors",
  "summary": "Researchers at CYFIRMA have identified a new ransomware strain named BARADAI that is actively targeting Windows systems. The malware encrypts files on local and network drives, appending a unique extension and dropping a ransom note with contact instructions. Observed victims are located in the USA, Brazil, France, Australia, Italy, Israel, and Malaysia, with a focus on the education, manufacturing, retail, and logistics sectors. Ransom demands have ranged from $10,000 to $80,000. BARADAI establishes persistence via `CurrentVersion\\Run` registry keys and manipulates system settings for defense evasion. Currently, no public decryption tool is available, and experts warn that future versions may adopt double extortion tactics.",
  "full_report": "## Executive Summary\nCybersecurity firm **[CYFIRMA](https://www.cyfirma.com/)** has reported the discovery of a new ransomware family dubbed **BARADAI**. This file-encrypting malware targets **[Windows](https://www.microsoft.com/en-us/windows)** systems, restricting access to data on local and network shares before demanding a ransom. The campaign appears to be global, with victims identified in the USA, Brazil, France, Australia, and other countries. Attackers are targeting a range of industries, including education, manufacturing, retail, and logistics, with ransom demands reported between $10,000 and $80,000. The malware ensures its persistence through registry modifications and employs defense evasion techniques. No free decryption tool is currently available, posing a significant data recovery challenge for victims.\n\n## Threat Overview\n**BARADAI ransomware** operates as a conventional file-encrypting threat. Upon execution, it systematically scans the victim's machine and any accessible network shares for files to encrypt. Once encrypted, files are appended with a unique extension, rendering them inaccessible. A ransom note is then created on the system, providing an email address for the victim to contact the operators and negotiate payment for a decryption key. The note includes a common threat that attempting to use third-party recovery tools will result in the permanent destruction of the data. The initial access vector has not been specified but likely involves common ransomware delivery methods such as phishing emails or exploitation of exposed remote services.\n\n## Technical Analysis\nBARADAI's functionality is straightforward but effective. Key technical aspects align with common ransomware TTPs:\n- **Execution & Impact:** The core of the attack is data encryption for impact ([`T1486 - Data Encrypted for Impact`](https://attack.mitre.org/techniques/T1486/)). It targets both local files and data on network shares ([`T1489 - Service Stop`](https://attack.mitre.org/techniques/T1489/) may be used to unlock files held by databases or applications).\n- **Persistence:** To survive a reboot, BARADAI creates an autorun entry in the `HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run` registry key, a classic persistence method ([`T1547.001 - Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder`](https://attack.mitre.org/techniques/T1547/001/)).\n- **Defense Evasion:** The report notes that the malware interacts with registry paths under `HKCU\\SOFTWARE` to store data and may manipulate internet settings and policy keys. This could be part of an effort to inhibit system recovery or evade detection ([`T1112 - Modify Registry`](https://attack.mitre.org/techniques/T1112/)).\n\nWhile not yet observed, it is highly probable that future versions will incorporate data exfiltration ([`T1041 - Exfiltration Over C2 Channel`](https://attack.mitre.org/techniques/T1041/)) before encryption to engage in double extortion, a tactic now standard for most ransomware groups.\n\n## Impact Assessment\nFor affected organizations, the primary impact is business disruption due to the unavailability of critical data and systems. The financial impact includes the cost of the ransom (if paid), recovery and remediation efforts, and lost revenue during downtime. The targeted sectors—manufacturing, retail, and logistics—are particularly vulnerable to operational halts. While currently not a data exfiltration threat, the potential for future versions to leak stolen data adds a layer of reputational and regulatory risk (e.g., GDPR, HIPAA) for victims.\n\n## IOCs — Directly from Articles\nNo specific technical Indicators of Compromise (IOCs) such as IP addresses, domains, or file hashes were mentioned in the source articles.\n\n## Cyber Observables — Hunting Hints\nSecurity teams can hunt for BARADAI and similar ransomware activity using these generic but effective patterns:\n| Type | Value | Description |\n|---|---|---|\n| `registry_key` | `HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run` | Monitor for new, unauthorized programs being added to this key for persistence. |\n| `file_name` | `*.baradai` (or other unique extension) | The appearance of a large number of files with a new, unknown extension is a strong indicator of ransomware. |\n| `file_name` | `ransom_note.txt` (or similar) | The creation of files with names like `readme.txt`, `decrypt_me.txt` in multiple directories is a hallmark of ransomware. |\n| `process_name` | High volume of file read/write/rename operations from a single process. | EDR solutions can detect the characteristic behavior of a process rapidly encrypting files across the file system. |\n\n## Detection & Response\n- **Detection:** Deploy EDR solutions with ransomware-specific behavioral detection capabilities. This includes monitoring for rapid file encryption activity (canary files), deletion of volume shadow copies (`vssadmin`), and modification of registry run keys. D3FEND's [`D3-FCR - File Content Rules`](https://d3fend.mitre.org/technique/d3f:FileContentRules) can be used to detect the creation of ransom notes.\n- **Response:** Upon detection of ransomware activity, immediately isolate the infected endpoints from the network to prevent lateral spread. If possible, power down affected systems to halt the encryption process, though this may complicate forensic analysis. Do not delete the encrypted files or the ransom note. Engage an incident response team to assess the scope and begin recovery from offline, immutable backups.\n\n## Mitigation\n- **Strategic:** The most effective mitigation against ransomware is a robust, multi-layered backup strategy ([`M1053 - Data Backup`](https://attack.mitre.org/mitigations/M1053/)). Follow the 3-2-1 rule: three copies of your data, on two different media types, with one copy stored offline and immutable. Implement a strong security awareness training program ([`M1017 - User Training`](https://attack.mitre.org/mitigations/M1017/)) to educate employees about phishing.\n- **Tactical:** Restrict user permissions to follow the principle of least privilege ([`M1026 - Privileged Account Management`](https://attack.mitre.org/mitigations/M1026/)). Ensure all software, especially on internet-facing systems like VPNs and RDP gateways, is fully patched ([`M1051 - Update Software`](https://attack.mitre.org/mitigations/M1051/)). Use email filtering and web gateways to block malicious attachments and links. Disable RDP on external-facing devices or secure it behind a VPN with MFA.",
  "twitter_post": "New ransomware strain 'BARADAI' discovered, targeting Windows systems in the US, Brazil, France & more. 💸 Demands range from $10k-$80k. Education, manufacturing, and retail sectors hit. No public decryptor available. #Ransomware #BARADAI #CyberCrime",
  "meta_description": "CYFIRMA researchers have identified a new ransomware strain called BARADAI, which targets Windows systems across multiple industries with ransom demands up to $80,000.",
  "category": [
    "Ransomware",
    "Malware"
  ],
  "severity": "high",
  "entities": [
    {
      "name": "CYFIRMA",
      "type": "security_organization",
      "url": "https://www.cyfirma.com/"
    },
    {
      "name": "BARADAI Ransomware",
      "type": "malware"
    },
    {
      "name": "Windows",
      "type": "product",
      "url": "https://www.microsoft.com/en-us/windows"
    }
  ],
  "cves": [],
  "sources": [
    {
      "url": "https://www.cyfirma.com/outofband/weekly-intelligence-report-08-may-2026/",
      "title": "Weekly Intelligence Report – 08 May 2026",
      "date": "2026-05-08",
      "friendly_name": "CYFIRMA",
      "website": "cyfirma.com"
    }
  ],
  "events": [],
  "mitre_techniques": [
    {
      "id": "T1486",
      "name": "Data Encrypted for Impact",
      "tactic": "Impact"
    },
    {
      "id": "T1547.001",
      "name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "tactic": "Persistence"
    },
    {
      "id": "T1112",
      "name": "Modify Registry",
      "tactic": "Defense Evasion"
    },
    {
      "id": "T1070.004",
      "name": "Indicator Removal: File Deletion",
      "tactic": "Defense Evasion"
    }
  ],
  "mitre_mitigations": [
    {
      "id": "M1053",
      "name": "Data Backup",
      "description": "Regularly backing up data to offline and immutable storage is the most effective defense, allowing for recovery without paying a ransom."
    },
    {
      "id": "M1017",
      "name": "User Training",
      "description": "Training users to recognize and report phishing attempts can prevent the initial access that often leads to ransomware infection."
    },
    {
      "id": "M1021",
      "name": "Restrict Web-Based Content",
      "description": "Using email and web filtering solutions to block malicious attachments, links, and scripts is a critical preventative measure."
    },
    {
      "id": "M1049",
      "name": "Antivirus/Antimalware",
      "description": "Modern EDR and antivirus solutions with behavioral analysis can detect and block ransomware activity based on its actions, such as rapid file encryption."
    }
  ],
  "d3fend_countermeasures": [
    {
      "technique_id": "D3-RB",
      "technique_name": "Remote Backup",
      "url": "https://d3fend.mitre.org/technique/d3f:RemoteBackup",
      "recommendation": "The ultimate defense against ransomware like BARADAI is ensuring data recoverability. Implementing a robust Remote Backup strategy is paramount. This involves more than just backing up data; it requires following the 3-2-1 rule: maintain at least three copies of data, on two different storage media, with at least one copy located off-site and offline or immutable. For BARADAI, which targets network shares, it's critical that backups are not accessible via the same credentials as the production network. Use separate authentication, pull-based backup mechanisms, and leverage cloud storage with object lock (immutability) features. Regularly test the restoration process to ensure backups are viable. This D3FEND technique neutralizes the primary leverage of encryption-only ransomware by making the ransom demand irrelevant.",
      "mitre_mitigation_id": "M1053"
    },
    {
      "technique_id": "D3-UAP",
      "technique_name": "User Account Permissions",
      "url": "https://d3fend.mitre.org/technique/d3f:UserAccountPermissions",
      "recommendation": "To limit the blast radius of a BARADAI infection, strict enforcement of User Account Permissions is essential. The principle of least privilege must be applied rigorously. Standard user accounts should not have administrative rights. More importantly, access to network shares should be tightly controlled. Users and service accounts should only have read/write access to the specific folders they absolutely need for their job function. Avoid using broad access groups like 'Everyone' or 'Authenticated Users'. By segmenting network file access, an infection on one user's machine will be contained, unable to encrypt data across the entire organization. This directly counters BARADAI's ability to traverse and encrypt 'accessible network shares'.",
      "mitre_mitigation_id": "M1018"
    },
    {
      "technique_id": "D3-FCR",
      "technique_name": "File Content Rules",
      "url": "https://d3fend.mitre.org/technique/d3f:FileContentRules",
      "recommendation": "A high-fidelity detection for ransomware like BARADAI can be achieved using File Content Rules, often deployed as 'honeyfiles' or 'canary files'. This involves strategically placing hidden files with specific, non-standard extensions (e.g., `.canary`) across file shares. An EDR or FIM solution is then configured to monitor these files for any modification, particularly rename operations. Since ransomware encrypts files indiscriminately, it will inevitably rename a canary file. This action should trigger an immediate, high-priority alert and, ideally, an automated response, such as isolating the source host from the network. This D3FEND technique provides a very early warning that an encryption process has begun, allowing for rapid response to minimize the overall damage.",
      "mitre_mitigation_id": "M1049"
    }
  ],
  "iocs": [],
  "cyber_observables": [
    {
      "type": "registry_key",
      "value": "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "description": "BARADAI creates entries in this registry key to establish persistence across reboots.",
      "context": "Registry monitoring, EDR, SIEM alerts",
      "confidence": "high"
    },
    {
      "type": "command_line_pattern",
      "value": "vssadmin.exe delete shadows /all /quiet",
      "description": "A common command used by ransomware to delete volume shadow copies and prevent easy system restore.",
      "context": "Windows Event ID 4688, EDR process monitoring",
      "confidence": "high"
    },
    {
      "type": "file_name",
      "value": "*.[unique_extension]",
      "description": "Widespread file renaming with a consistent, unknown extension is a primary indicator of a ransomware attack in progress.",
      "context": "File integrity monitoring (FIM), EDR",
      "confidence": "high"
    },
    {
      "type": "log_source",
      "value": "Windows Security Event Log",
      "description": "Monitor for a high rate of file modification events (Event ID 4663) originating from a single, non-system process.",
      "context": "SIEM correlation rules",
      "confidence": "medium"
    }
  ],
  "tags": [
    "Ransomware",
    "BARADAI",
    "CYFIRMA",
    "Malware",
    "Windows",
    "Data Encryption"
  ],
  "extract_datetime": "2026-05-08T15:00:00.000Z",
  "article_type": "NewsArticle",
  "impact_scope": {
    "geographic_scope": "global",
    "countries_affected": [
      "USA",
      "Brazil",
      "France",
      "Australia",
      "Italy",
      "Israel",
      "Malaysia"
    ],
    "industries_affected": [
      "Education",
      "Manufacturing",
      "Retail",
      "Other"
    ]
  },
  "pub_date": "2026-05-08",
  "reading_time_minutes": 5,
  "createdAt": "2026-05-08T15:00:00.000Z",
  "updatedAt": "2026-05-08T15:00:00.000Z"
}