# data/articles/stealthy-python-backdoor-deepdoor-steals-credentials-via-tunneling.json
{
  "id": "1679065f-32df-4bf0-bfe2-36ba0f5a7b9a",
  "slug": "stealthy-python-backdoor-deepdoor-steals-credentials-via-tunneling",
  "headline": "New Stealthy Python Backdoor 'DEEP#DOOR' Steals Credentials Using Legitimate Tunneling Service",
  "title": "Researchers Uncover 'DEEP#DOOR', a Stealthy Python Backdoor Stealing Credentials",
  "summary": "Security researchers from Securonix have discovered a sophisticated, multi-stage Python backdoor named DEEP#DOOR. The malware operates as a full-featured Remote Access Trojan (RAT), beginning with an obfuscated batch script that embeds the Python payload. It disables Windows security controls, establishes persistence through multiple methods (Startup folder, Run keys, Scheduled Tasks, WMI), and uses the legitimate 'bore.pub' tunneling service for stealthy command-and-control and data exfiltration, focusing on browser and cloud credentials.",
  "full_report": "## Executive Summary\n\nResearchers at **[Securonix](https://www.securonix.com/)** have detailed a new and stealthy Python-based backdoor framework they have named **DEEP#DOOR**. This malware is a full-featured Remote Access Trojan (RAT) designed to gain persistent access to Windows systems and exfiltrate sensitive data. The attack chain begins with an obfuscated batch script that contains the core Python implant, a technique that minimizes its network footprint. Upon execution, **DEEP#DOOR** disables Windows security features and establishes persistence using at least four different methods to ensure its survival. A key feature of the malware is its use of `bore.pub`, a legitimate tunneling service, for its command and control (C2) communications. This allows it to bypass network security controls and stealthily exfiltrate stolen browser and cloud credentials.\n\n---\n\n## Threat Overview\n\n**DEEP#DOOR** is a sophisticated RAT that demonstrates a trend towards fileless and script-based malware that abuses legitimate system tools and services to evade detection. The initial access vector is believed to be phishing or a similar social engineering tactic to trick a user into running the initial dropper.\n\nThe malware's primary objectives are to:\n1.  Establish a persistent foothold on the victim's machine.\n2.  Disable security software.\n3.  Steal credentials, particularly from web browsers and cloud service applications.\n4.  Provide the attacker with remote access for further post-exploitation activities.\n\n## Technical Analysis\n\nThe attack chain unfolds in several stages:\n\n1.  **Initial Dropper (`install_obf.bat`):** The attack starts with an obfuscated batch script. This script contains the embedded and encoded Python payload, avoiding the need to download a second stage from the internet.\n2.  **Defense Evasion:** The batch script first attempts to disable Windows security controls, likely including Windows Defender and other monitoring services.\n3.  **Payload Extraction:** The script then decodes and writes the embedded Python implant (`svc.py`) to the disk and executes it.\n4.  **Persistence ([`T1547.001`](https://attack.mitre.org/techniques/T1547/001/), [`T1543.003`](https://attack.mitre.org/techniques/T1543/003/)):** **DEEP#DOOR** is highly resilient, establishing persistence through multiple redundant methods simultaneously:\n    - Creating scripts in the user's Startup folder.\n    - Adding entries to the `Registry Run` keys.\n    - Creating a `Scheduled Task`.\n    - Using `Windows Management Instrumentation (WMI)` event subscriptions.\n5.  **Command and Control ([`T1572`](https://attack.mitre.org/techniques/T1572/)):** The Python implant communicates with `bore.pub`, a legitimate, Rust-based tunneling service. It uses this service to create a reverse tunnel from the victim's machine to the attacker's C2 server. This makes the C2 traffic appear as legitimate outbound web traffic to a known service, helping it bypass firewalls and network monitoring tools.\n6.  **Credential Theft ([`T1555`](https://attack.mitre.org/techniques/T1555/)):** Once active, the RAT can receive commands to execute various modules, including those designed to steal credentials from browsers, email clients, and cloud applications.\n\n## Impact Assessment\n\nA successful **DEEP#DOOR** infection can have severe consequences:\n- **Complete System Compromise:** As a RAT, it gives the attacker full control over the infected machine, allowing them to install other malware (like ransomware), spy on the user, and use the machine as a pivot point for lateral movement.\n- **Widespread Credential Theft:** The theft of saved browser and cloud credentials can lead to the compromise of numerous other online accounts, both personal and corporate.\n- **Data Exfiltration:** The attacker can use their remote access to search for and exfiltrate any sensitive files or data on the compromised system or accessible network shares.\n- **Difficult Detection:** The use of legitimate services (`bore.pub`), script-based execution (`Python`, `Batch`), and multiple persistence mechanisms makes the malware difficult to detect and fully eradicate.\n\n---\n\n## IOCs — Directly from Articles\n\n| Type | Value | Description |\n|---|---|---|\n| File Name | `install_obf.bat` | The initial obfuscated batch script dropper. |\n| File Name | `svc.py` | The name of the core Python implant payload. |\n| Domain | `bore.pub` | The legitimate tunneling service abused for C2 communications. |\n\n## Detection & Response\n\n1.  **Monitor for Tunneling Services:** Monitor network traffic for connections to known tunneling services like `bore.pub`, `ngrok.io`, or `localtunnel.me`. While these have legitimate uses, an unexpected connection from a user workstation or server is highly suspicious and should be investigated.\n2.  **Script-Block Logging:** Enable PowerShell and other script-block logging and forward the logs to a SIEM. This can capture the content of malicious scripts like the initial batch file, even if they are obfuscated.\n3.  **Persistence Mechanism Auditing ([D3-SFA](https://d3fend.mitre.org/technique/d3f:SystemFileAnalysis)):** Regularly audit the common persistence locations used by **DEEP#DOOR** (Startup folders, Run keys, Scheduled Tasks, WMI subscriptions) for any unauthorized or suspicious entries.\n\n## Mitigation\n\n1.  **Application Control ([D3-EAL](https://d3fend.mitre.org/technique/d3f:ExecutableAllowlisting)):** Use application control solutions like AppLocker to restrict the execution of batch scripts and unsigned Python scripts in user-writable directories. This can prevent the initial dropper from running.\n2.  **Endpoint Protection:** Deploy a modern EDR solution capable of detecting and blocking suspicious script execution and fileless malware techniques.\n3.  **Egress Traffic Filtering ([D3-OTF](https://d3fend.mitre.org/technique/d3f:OutboundTrafficFiltering)):** If your organization does not have a business need for tunneling services, explicitly block the domains of services like `bore.pub` at your web proxy or firewall.\n4.  **User Training:** Educate users about the dangers of running unsolicited scripts or attachments from emails.",
  "twitter_post": "🐍 New Python backdoor 'DEEP#DOOR' discovered! This stealthy RAT disables security, uses multiple persistence methods, and exfiltrates credentials via the legit 'bore.pub' tunneling service. #Malware #Python #Backdoor #ThreatIntel",
  "meta_description": "Securonix researchers detail DEEP#DOOR, a new stealthy Python backdoor that functions as a RAT to steal credentials from Windows systems by abusing a legitimate tunneling service for C2.",
  "category": [
    "Malware",
    "Threat Intelligence",
    "Cyberattack"
  ],
  "severity": "high",
  "entities": [
    {
      "name": "DEEP#DOOR",
      "type": "malware"
    },
    {
      "name": "Securonix",
      "type": "security_organization",
      "url": "https://www.securonix.com/"
    },
    {
      "name": "Python",
      "type": "technology"
    },
    {
      "name": "Windows",
      "type": "product"
    }
  ],
  "cves": [],
  "sources": [
    {
      "url": "https://thehackernews.com/2026/04/new-python-backdoor-uses-tunneling.html",
      "title": "New Python Backdoor Uses Tunneling Service to Steal Browser and Cloud Credentials",
      "date": "2026-04-30",
      "friendly_name": "The Hacker News"
    },
    {
      "url": "https://www.cimetrics.com/blog-posts/2026/04/29/april-2026-cybersecurity-news",
      "title": "April 2026 Cybersecurity News",
      "date": "2026-04-29",
      "friendly_name": "Cimetrics"
    }
  ],
  "events": [],
  "mitre_techniques": [
    {
      "id": "T1059.006",
      "name": "Python",
      "tactic": "Execution"
    },
    {
      "id": "T1572",
      "name": "Protocol Tunneling",
      "tactic": "Command and Control"
    },
    {
      "id": "T1547.001",
      "name": "Registry Run Keys / Startup Folder",
      "tactic": "Persistence"
    },
    {
      "id": "T1546.003",
      "name": "WMI Event Subscription",
      "tactic": "Persistence"
    },
    {
      "id": "T1555",
      "name": "Credentials from Password Stores",
      "tactic": "Credential Access"
    },
    {
      "id": "T1562.001",
      "name": "Disable or Modify Tools",
      "tactic": "Defense Evasion"
    }
  ],
  "mitre_mitigations": [
    {
      "id": "M1038",
      "name": "Execution Prevention",
      "description": "Using application control policies to restrict the execution of scripting languages like Python and Batch from untrusted locations can prevent the malware from running.",
      "domain": "enterprise"
    },
    {
      "id": "M1037",
      "name": "Filter Network Traffic",
      "description": "Explicitly blocking outbound connections to known tunneling services like bore.pub at the network perimeter can disrupt the malware's C2 channel.",
      "domain": "enterprise"
    },
    {
      "id": "M1047",
      "name": "Audit",
      "description": "Regularly auditing common persistence locations (Registry, Startup folders, WMI) can help detect the presence of DEEP#DOOR and similar threats.",
      "domain": "enterprise"
    }
  ],
  "d3fend_countermeasures": [
    {
      "technique_id": "D3-OTF",
      "technique_name": "Outbound Traffic Filtering",
      "url": "https://d3fend.mitre.org/technique/d3f:OutboundTrafficFiltering",
      "recommendation": "A key feature of DEEP#DOOR is its abuse of the legitimate service 'bore.pub' for C2 tunneling. Organizations should implement an egress filtering policy that blocks access to known tunneling services unless there is an explicit business justification. Create a category of 'Proxy and Tunneling Services' in your web proxy or firewall and set it to block. This list should include `bore.pub`, `ngrok.io`, and other similar services. This simple network control directly severs the malware's communication channel, rendering it unable to receive commands or exfiltrate data.",
      "mitre_mitigation_id": "M1037"
    },
    {
      "technique_id": "D3-EAL",
      "technique_name": "Executable Allowlisting",
      "url": "https://d3fend.mitre.org/technique/d3f:ExecutableAllowlisting",
      "recommendation": "To prevent the initial execution of the DEEP#DOOR dropper, organizations should implement application allowlisting using tools like Windows AppLocker. Configure a policy that prevents the execution of batch files (`.bat`) and Python scripts (`.py`) from user-writable directories such as `C:\\Users\\` and `%APPDATA%`. By only allowing scripts to run from protected, administrator-controlled directories (e.g., `C:\\Program Files\\`), you can effectively block this and many other script-based threats that rely on a user downloading and running a file in their local profile.",
      "mitre_mitigation_id": "M1038"
    }
  ],
  "iocs": [
    {
      "type": "file_name",
      "value": "install_obf.bat",
      "description": "Initial dropper script"
    },
    {
      "type": "file_name",
      "value": "svc.py",
      "description": "Core Python payload"
    },
    {
      "type": "domain",
      "value": "bore.pub",
      "description": "Legitimate tunneling service used for C2"
    }
  ],
  "cyber_observables": [
    {
      "type": "domain",
      "value": "bore.pub",
      "description": "Outbound connections to this legitimate tunneling service from endpoints that have no business reason to use it is a strong indicator of compromise.",
      "context": "DNS logs, proxy logs, firewall logs.",
      "confidence": "high"
    },
    {
      "type": "process_name",
      "value": "python.exe",
      "description": "Execution of the Python interpreter from an unusual file path or by an unexpected parent process (e.g., cmd.exe running a batch script).",
      "context": "EDR process monitoring, Windows Event ID 4688.",
      "confidence": "medium"
    },
    {
      "type": "registry_key",
      "value": "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "description": "DEEP#DOOR adds entries to this registry key to maintain persistence across reboots. Monitor for new, suspicious entries.",
      "context": "Registry monitoring tools, EDR, or SIEM.",
      "confidence": "high"
    },
    {
      "type": "service_name",
      "value": "WMI Event Subscription",
      "description": "The malware uses WMI event subscriptions for fileless persistence. Hunting for newly created WMI consumers, filters, and bindings is a key detection strategy.",
      "context": "PowerShell scripts for querying WMI, specialized EDR WMI monitoring.",
      "confidence": "high"
    }
  ],
  "tags": [
    "DEEP#DOOR",
    "Python",
    "Backdoor",
    "RAT",
    "Securonix",
    "WMI Persistence",
    "Tunneling"
  ],
  "extract_datetime": "2026-04-30T15:00:00.000Z",
  "impact_scope": {
    "geographic_scope": "global",
    "industries_affected": [
      "Other"
    ],
    "other_affected": [
      "Windows system users"
    ]
  },
  "pub_date": "2026-04-30",
  "reading_time_minutes": 5,
  "createdAt": "2026-04-30T15:00:00.000Z",
  "updatedAt": "2026-04-30T15:00:00.000Z"
}