# stix/china-apt-mirrorface-expands-targeting-to-eu-diplomats-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T14:38:17.670Z",
    "article_slug": "china-apt-mirrorface-expands-targeting-to-eu-diplomats",
    "article_date": "2026-04-27",
    "token_usage": {
      "input_tokens": 137536,
      "output_tokens": 2562,
      "cached_tokens": 132145,
      "total_tokens": 140098,
      "cache_hit": true,
      "estimated_cost_usd": 0.048877
    }
  },
  "report_name": "China-Linked APT MirrorFace Expands Targeting to EU Diplomatic Entity",
  "threat_actors": [
    {
      "name": "MirrorFace",
      "mitre_id": "",
      "aliases": [],
      "description": "A China-aligned advanced persistent threat (APT) group focused on cyber-espionage. The group has historically targeted entities in Japan but was observed expanding its operations to a diplomatic organization in the European Union in April 2024. MirrorFace is known for using custom malware and legitimate tools like SoftEther VPN for persistence and C2.",
      "goals": "Intelligence gathering related to political and economic affairs.",
      "primary_motivation": "organizational-gain",
      "sophistication": "advanced"
    }
  ],
  "attack_patterns": [
    {
      "name": "Protocol Tunneling",
      "mitre_id": "T1572",
      "description": "The threat actor used SoftEther VPN to tunnel command and control traffic, making it difficult to distinguish from normal network activity."
    },
    {
      "name": "Application Layer Protocol: Web Protocols",
      "mitre_id": "T1071.001",
      "description": "Command and control traffic, tunneled via SoftEther VPN, was disguised as normal HTTPS traffic."
    }
  ],
  "malware": [],
  "tools": [
    {
      "name": "SoftEther VPN",
      "mitre_id": "",
      "aliases": [],
      "description": "Leveraged by MirrorFace and other China-aligned groups to maintain stealthy, persistent access to compromised networks and for command and control (C2) communications.",
      "tool_types": [
        "remote-access"
      ]
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "European Union",
      "identity_class": "organization",
      "sectors": [
        "government"
      ],
      "description": "A diplomatic organization within the European Union was targeted by the MirrorFace APT group, marking a shift from the group's traditional focus on Japan."
    }
  ],
  "locations": [
    {
      "name": "China",
      "region": "eastern-asia"
    },
    {
      "name": "Japan",
      "region": "eastern-asia"
    },
    {
      "name": "European Union",
      "region": "europe"
    }
  ],
  "infrastructure": [
    {
      "name": "Unknown VPN Servers",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "SoftEther VPN was used to establish persistent connections to non-corporate VPN endpoints for command and control."
    }
  ],
  "observables": [
    {
      "type": "file",
      "name": "sevpn.exe"
    },
    {
      "type": "file",
      "name": "sevpnclient.exe"
    }
  ],
  "indicators": [
    {
      "name": "SoftEther VPN Client Execution",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[file:name IN ('sevpn.exe', 'sevpnclient.exe')]",
      "x_netsecio_sigma_rule": "title: Unauthorized SoftEther VPN Client Execution\nid: 2b1a9c8d-4e5f-4a3b-8c7d-6e5f4a3b2c1d\nstatus: experimental\ndescription: Detects the execution of SoftEther VPN client processes, which could indicate unauthorized remote access by threat actors like MirrorFace.\nreferences:\n  - https://www.eset.com/int/about/newsroom/press-releases/research/eset-releases-latest-apt-report-china-aligned-groups-expand-targeting-iran-advances-diplomatic-esp/\nauthor: cyber.netsecops.io\ndate: 2026-04-27\ntags:\n  - attack.t1219\n  - attack.t1572\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith:\n      - '\\sevpn.exe'\n      - '\\sevpnclient.exe'\n  condition: selection\nfalsepositives:\n  - Legitimate use of SoftEther VPN.\n",
      "pattern": "",
      "description": "Detects the execution of SoftEther VPN client processes. The presence of these executables on workstations or servers where they are not authorized is a strong indicator of compromise, associated with groups like MirrorFace."
    },
    {
      "name": "Suspicious LNK File in Startup Folder",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious LNK File Creation in Startup Folder\nid: 9f8e7d6c-5b4a-4987-b6a5-4321fedcba98\nstatus: experimental\ndescription: Detects the creation of LNK files in user or system startup directories, a common persistence technique used by APTs.\nreferences:\n  - https://www.eset.com/int/about/newsroom/press-releases/research/eset-releases-latest-apt-report-china-aligned-groups-expand-targeting-iran-advances-diplomatic-esp/\nauthor: cyber.netsecops.io\ndate: 2026-04-27\ntags:\n  - attack.t1547.009\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: file_event\n  product: windows\ndetection:\n  selection:\n    TargetFilename|endswith: '.lnk'\n    TargetFilename|contains:\n      - '\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\'\n      - '\\ProgramData\\Microsoft\\Windows\\Start Menu\\Programs\\StartUp\\'\n  condition: selection\nfalsepositives:\n  - Legitimate software installers creating startup shortcuts.\n",
      "description": "Detects the creation of LNK files in startup directories, a common persistence technique suggested for hunting MirrorFace activity."
    },
    {
      "name": "Suspicious Registry Run Key Creation",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Registry Run Key Creation\nid: 1a2b3c4d-5e6f-7a8b-9c0d-1e2f3a4b5c6d\nstatus: experimental\ndescription: Detects modification of the registry run keys, a common persistence technique used by malware and APTs.\nreferences:\n  - https://www.eset.com/int/about/newsroom/press-releases/research/eset-releases-latest-apt-report-china-aligned-groups-expand-targeting-iran-advances-diplomatic-esp/\nauthor: cyber.netsecops.io\ndate: 2026-04-27\ntags:\n  - attack.t1547.001\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains:\n      - 'Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n      - 'Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce'\n    EventType: 'SetValue'\n  condition: selection\nfalsepositives:\n  - Legitimate software installers and updaters.\n",
      "description": "Detects the creation or modification of registry run keys for persistence, a hunting hint provided for detecting APT activity like MirrorFace."
    }
  ],
  "relationships": [
    {
      "source_name": "MirrorFace",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "SoftEther VPN",
      "target_type": "tool"
    },
    {
      "source_name": "MirrorFace",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "European Union",
      "target_type": "identity"
    },
    {
      "source_name": "MirrorFace",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Japan",
      "target_type": "location"
    },
    {
      "source_name": "MirrorFace",
      "source_type": "threat-actor",
      "relationship_type": "originates-from",
      "target_name": "China",
      "target_type": "location"
    },
    {
      "source_name": "SoftEther VPN",
      "source_type": "tool",
      "relationship_type": "uses",
      "target_name": "Protocol Tunneling",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SoftEther VPN",
      "source_type": "tool",
      "relationship_type": "uses",
      "target_name": "Application Layer Protocol: Web Protocols",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SoftEther VPN",
      "source_type": "tool",
      "relationship_type": "communicates-with",
      "target_name": "Unknown VPN Servers",
      "target_type": "infrastructure"
    },
    {
      "source_name": "SoftEther VPN Client Execution",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "sevpn.exe",
      "target_type": "file"
    },
    {
      "source_name": "SoftEther VPN Client Execution",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "sevpnclient.exe",
      "target_type": "file"
    },
    {
      "source_name": "SoftEther VPN Client Execution",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "MirrorFace",
      "target_type": "threat-actor"
    }
  ]
}