# stix/chinese-apt-mustang-panda-targets-indian-banks-and-korean-policy-experts-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T17:13:01.714Z",
    "article_slug": "chinese-apt-mustang-panda-targets-indian-banks-and-korean-policy-experts",
    "article_date": "2026-04-21",
    "token_usage": {
      "input_tokens": 141943,
      "output_tokens": 3126,
      "cached_tokens": 136033,
      "total_tokens": 145069,
      "cache_hit": true,
      "estimated_cost_usd": 0.055652
    }
  },
  "report_name": "Chinese APT Mustang Panda Targets Indian Banks, Korean Policy Experts in Espionage Campaign",
  "threat_actors": [
    {
      "name": "Mustang Panda",
      "mitre_id": "G0129",
      "aliases": [
        "TA416",
        "Bronze President",
        "Stately Taurus"
      ],
      "description": "A Chinese state-sponsored cyber-espionage group conducting a campaign targeting the Indian financial sector and public policy experts in Korea and the U.S. The group uses spear-phishing and DLL side-loading to deploy the LotusLite backdoor for intelligence gathering.",
      "goals": "Intelligence gathering aligned with Beijing's geopolitical interests.",
      "primary_motivation": "organizational-gain",
      "resource_level": "government",
      "sophistication": "intermediate"
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Attachment",
      "mitre_id": "T1566.001",
      "description": "Used as the initial access vector, delivering malicious files like ZIP archives with LNK files to victims via email."
    },
    {
      "name": "Hijack Execution Flow",
      "mitre_id": "T1574.002",
      "description": "The core execution technique used by Mustang Panda, where a legitimate, signed application is tricked into loading a malicious DLL placed in the same directory."
    },
    {
      "name": "Boot or Logon Autostart Execution",
      "mitre_id": "T1547.001",
      "description": "The LotusLite backdoor establishes persistence by creating or modifying a Windows Registry key in a 'Run' key to ensure it executes on system startup."
    },
    {
      "name": "Application Layer Protocol",
      "mitre_id": "T1071.001",
      "description": "The LotusLite backdoor uses web protocols to connect to its command and control server for remote command execution and file transfer."
    },
    {
      "name": "Windows Command Shell",
      "mitre_id": "T1059.003",
      "description": "The LotusLite backdoor provides the attacker with a remote shell to execute commands on the compromised system."
    },
    {
      "name": "Obfuscated Files or Information",
      "mitre_id": "T1027",
      "description": "The malware used decoys, such as displaying a pop-up with 'HDFC Bank' branding, to hide its true purpose and allay victim suspicion."
    }
  ],
  "malware": [
    {
      "name": "LotusLite",
      "mitre_id": "",
      "type": "backdoor",
      "is_family": false,
      "aliases": [],
      "description": "A custom backdoor deployed by Mustang Panda via DLL side-loading. It enables remote command execution, file access, and intelligence gathering from compromised systems.",
      "capabilities": [
        "communicates-with-c2",
        "controls-local-machine",
        "persists-after-system-reboot"
      ]
    }
  ],
  "tools": [],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "HDFC Bank",
      "identity_class": "organization",
      "sectors": [
        "financial-services"
      ],
      "description": "An Indian private-sector bank whose branding was used as a decoy by Mustang Panda's malware to deceive victims in the Indian financial sector."
    },
    {
      "name": "Victor Cha",
      "identity_class": "individual",
      "sectors": [
        "education",
        "government"
      ],
      "description": "An American political scientist who was impersonated by Mustang Panda using a Google account to add legitimacy to spear-phishing outreach targeting policy experts."
    }
  ],
  "locations": [
    {
      "name": "China",
      "region": "eastern-asia"
    },
    {
      "name": "India",
      "region": "southern-asia"
    },
    {
      "name": "Korea",
      "region": "eastern-asia"
    },
    {
      "name": "United States",
      "region": "northern-america"
    }
  ],
  "infrastructure": [
    {
      "name": "LotusLite C2 Server",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "Attacker-controlled command and control server used by the LotusLite backdoor for remote command execution and file access."
    }
  ],
  "observables": [
    {
      "type": "windows-registry-key",
      "value": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
    }
  ],
  "indicators": [
    {
      "name": "Mustang Panda DLL Side-Loading Detection",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Potential DLL Side-Loading by Legitimate Process\nid: 2a9b3c4d-5e6f-4a7b-8c9d-0e1f2a3b4c5d\nstatus: experimental\ndescription: Detects a legitimate, signed process loading a DLL from a non-standard, user-writable directory, a technique commonly used by Mustang Panda.\nreferences:\n  - https://www.acronis.com/en-us/blog/posts/mustang-panda-attack-chain\nauthor: cyber.netsecops.io\ndate: 2026-04-21\ntags:\n  - attack.t1574.002\n  - attack.persistence\n  - attack.defense_evasion\nlevel: medium\nlogsource:\n  category: image_load\n  product: windows\ndetection:\n  selection_process:\n    Image|endswith:\n      - '\\Acrobat.exe'\n      - '\\WINWORD.EXE'\n      - '\\EXCEL.EXE'\n      - '\\POWERPNT.EXE'\n  selection_path:\n    ImageLoaded|startswith:\n      - 'C:\\Users\\'\n      - 'C:\\ProgramData\\'\n      - 'C:\\Temp\\'\n      - 'C:\\Windows\\Temp\\'\n  condition: all of selection_*\nfalsepositives:\n  - Legitimate applications that incorrectly package their DLLs.\n  - Portable applications running from user directories.",
      "description": "Detects the DLL side-loading technique used by Mustang Panda where a legitimate application loads a malicious DLL from a non-standard directory."
    },
    {
      "name": "Mustang Panda Registry Run Key Persistence",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[windows-registry-key:key = 'HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run']",
      "x_netsecio_sigma_rule": "title: Suspicious Modification of HKCU Run Key\nid: 1b8c2d3e-4f5a-6b7c-8d9e-0f1a2b3c4d5e\nstatus: experimental\ndescription: Detects modification of the HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run registry key, a common persistence technique used by malware like LotusLite.\nreferences:\n  - https://www.acronis.com/en-us/blog/posts/mustang-panda-attack-chain\nauthor: cyber.netsecops.io\ndate: 2026-04-21\ntags:\n  - attack.t1547.001\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    EventType: 'SetValue'\n    TargetObject|contains: '\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n  condition: selection\nfalsepositives:\n  - Legitimate software installers adding autostart entries.",
      "pattern": "",
      "description": "Detects the creation of persistence mechanisms in the Windows Registry Run key, as used by LotusLite."
    }
  ],
  "relationships": [
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "LotusLite",
      "target_type": "malware"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Spearphishing Attachment",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Hijack Execution Flow",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "India",
      "target_type": "location"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Korea",
      "target_type": "location"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "United States",
      "target_type": "location"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "impersonates",
      "target_name": "Victor Cha",
      "target_type": "identity"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "originates-from",
      "target_name": "China",
      "target_type": "location"
    },
    {
      "source_name": "LotusLite",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Boot or Logon Autostart Execution",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "LotusLite",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Windows Command Shell",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "LotusLite",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Application Layer Protocol",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "LotusLite",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Obfuscated Files or Information",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "LotusLite",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "LotusLite C2 Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Mustang Panda",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "HDFC Bank",
      "target_type": "identity"
    },
    {
      "source_name": "Mustang Panda DLL Side-Loading Detection",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Hijack Execution Flow",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mustang Panda Registry Run Key Persistence",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Boot or Logon Autostart Execution",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mustang Panda Registry Run Key Persistence",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "target_type": "windows-registry-key"
    }
  ]
}