# stix/fsb-linked-gamaredon-apt-deploys-fileless-worm-against-ukraine-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-06-01T16:06:20.616Z",
    "article_slug": "fsb-linked-gamaredon-apt-deploys-fileless-worm-against-ukraine",
    "article_date": "2026-06-01",
    "token_usage": {
      "input_tokens": 143403,
      "output_tokens": 3800,
      "cached_tokens": 137782,
      "total_tokens": 147203,
      "cache_hit": true,
      "estimated_cost_usd": 0.062249
    }
  },
  "report_name": "FSB-Linked Gamaredon APT Deploys Stealthy Fileless Worm Against Ukrainian Targets",
  "threat_actors": [
    {
      "name": "Gamaredon",
      "mitre_id": "G0047",
      "aliases": [
        "Shuckworm",
        "Primitive Bear"
      ],
      "description": "A Russian state-sponsored threat actor, attributed to the FSB, that deployed a new fileless VBScript worm in a cyber-espionage campaign against Ukraine. The group focuses on long-term intelligence gathering and data theft from high-value targets.",
      "goals": "Long-term intelligence gathering and data theft from Ukrainian government, military, and critical infrastructure entities.",
      "primary_motivation": "organizational-gain",
      "resource_level": "government",
      "sophistication": "advanced"
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Attachment",
      "mitre_id": "T1566.001",
      "description": "Gamaredon uses spear-phishing emails with malicious attachments, such as LNK files or documents with remote templates, for initial access."
    },
    {
      "name": "NTFS File Attributes",
      "mitre_id": "T1564.004",
      "description": "The fileless VBScript worm hides its malicious code within an NTFS Alternate Data Stream (ADS) of a legitimate-looking file to evade detection."
    },
    {
      "name": "Windows Management Instrumentation",
      "mitre_id": "T1047",
      "description": "The malware uses WMI event subscriptions to execute the VBScript payload directly from its hidden NTFS alternate data stream."
    },
    {
      "name": "Scheduled Task",
      "mitre_id": "T1053.005",
      "description": "The malware uses scheduled tasks to execute the VBScript payload directly from its hidden NTFS alternate data stream for persistence and execution."
    },
    {
      "name": "Visual Basic",
      "mitre_id": "T1059.005",
      "description": "The fileless worm is primarily written in VBScript, which is executed by the Windows Script Host (wscript.exe)."
    },
    {
      "name": "Automated Collection",
      "mitre_id": "T1119",
      "description": "The worm's payload is designed to search for documents with specific keywords or file types on compromised systems."
    },
    {
      "name": "Exfiltration Over C2 Channel",
      "mitre_id": "T1041",
      "description": "Collected documents and data are exfiltrated to an attacker-controlled command and control server."
    }
  ],
  "malware": [
    {
      "name": "VBScript Worm",
      "mitre_id": "",
      "type": "worm",
      "is_family": false,
      "aliases": [],
      "description": "A new, fileless VBScript worm deployed by Gamaredon that uses NTFS alternate data streams (ADS) for stealth and persistence. Its primary purpose is espionage and data exfiltration.",
      "capabilities": [
        "hides-artifacts",
        "persists-after-system-reboot",
        "exfiltrates-data",
        "communicates-with-c2"
      ]
    }
  ],
  "tools": [
    {
      "name": "wscript.exe",
      "mitre_id": "",
      "aliases": [],
      "description": "The Windows Script Host is used to execute the malicious VBScript payload stored in an NTFS alternate data stream.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "Ukrainian government organizations",
      "identity_class": "class",
      "sectors": [
        "government"
      ],
      "description": "A class of organizations in Ukraine, including government bodies, targeted by the Gamaredon APT group for cyber-espionage."
    },
    {
      "name": "Ukrainian military",
      "identity_class": "class",
      "sectors": [
        "defense"
      ],
      "description": "The Ukrainian military, targeted by the Gamaredon APT group for long-term intelligence gathering."
    },
    {
      "name": "Ukrainian critical infrastructure",
      "identity_class": "class",
      "sectors": [
        "infrastructure"
      ],
      "description": "Critical infrastructure entities in Ukraine targeted by the Gamaredon APT group in an ongoing cyber-espionage campaign."
    },
    {
      "name": "Federal Security Service (FSB)",
      "identity_class": "organization",
      "sectors": [
        "government"
      ],
      "description": "The Russian Federal Security Service, to which the Gamaredon APT group is attributed."
    }
  ],
  "locations": [
    {
      "name": "Russia",
      "country": "Russia",
      "region": "eastern-europe"
    },
    {
      "name": "Ukraine",
      "country": "Ukraine",
      "region": "eastern-europe"
    }
  ],
  "infrastructure": [],
  "observables": [],
  "indicators": [
    {
      "name": "Execution from NTFS Alternate Data Stream",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Execution from NTFS Alternate Data Stream\nid: 2b1a9c8d-4e5f-4a3b-8c7d-6e5f4a3b2c1d\nstatus: experimental\ndescription: Detects execution of scripts from NTFS alternate data streams, a technique used by Gamaredon's VBScript worm.\nreferences:\n  - https://www.infosecurity-magazine.com/news/fsb-gamaredon-hides-worm-windows/\nauthor: cyber.netsecops.io\ndate: 2026-06-01\ntags:\n  - attack.t1564.004\n  - attack.t1059.005\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith:\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n    CommandLine|contains: ':'\n  filter_legit:\n    CommandLine|re: '^[a-zA-Z]:\\\\'\n  condition: selection and not filter_legit",
      "description": "Detects the execution of a script host where the command line contains a colon in the path, indicative of executing from an NTFS alternate data stream. This was a key TTP of the Gamaredon worm."
    },
    {
      "name": "Enumeration of NTFS Alternate Data Streams via Dir",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Enumeration of NTFS Alternate Data Streams via Dir\nid: f8c7b6a5-4d3c-4b2a-9a8b-7c6d5e4f3a2b\nstatus: experimental\ndescription: 'Detects the use of ''dir /r'' to discover NTFS alternate data streams. This may be used by adversaries for reconnaissance or by defenders for hunting.'\nreferences:\n  - https://www.infosecurity-magazine.com/news/fsb-gamaredon-hides-worm-windows/\nauthor: cyber.netsecops.io\ndate: 2026-06-01\ntags:\n  - attack.t1564.004\nlevel: low\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\dir.exe'\n    CommandLine|contains: ' /r'\n  condition: selection",
      "description": "Detects the use of the 'dir /r' command, which is used to reveal NTFS alternate data streams. This is a hunting hint provided in the article to find the hidden malware components."
    },
    {
      "name": "Suspicious Parent for Windows Script Host",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Parent for Windows Script Host\nid: 1a2b3c4d-5e6f-4a9b-8c7d-6e5f4a3b2c1d\nstatus: experimental\ndescription: Monitors for wscript.exe or cscript.exe being launched by suspicious parent processes like svchost.exe (for scheduled tasks) or WmiPrvSE.exe (for WMI), a TTP of Gamaredon.\nreferences:\n  - https://www.infosecurity-magazine.com/news/fsb-gamaredon-hides-worm-windows/\nauthor: cyber.netsecops.io\ndate: 2026-06-01\ntags:\n  - attack.t1059.005\n  - attack.t1047\n  - attack.t1053.005\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_child:\n    Image|endswith:\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n  selection_parent:\n    ParentImage|endswith:\n      - '\\svchost.exe'\n      - '\\WmiPrvSE.exe'\n  condition: all of selection_*",
      "description": "Detects Windows Script Host processes being launched by parent processes associated with scheduled tasks (svchost.exe) or WMI (WmiPrvSE.exe), as described in the article's analysis of the malware's execution."
    }
  ],
  "relationships": [
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "attributed-to",
      "target_name": "Federal Security Service (FSB)",
      "target_type": "identity"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "VBScript Worm",
      "target_type": "malware"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Spearphishing Attachment",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukrainian government organizations",
      "target_type": "identity"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukrainian military",
      "target_type": "identity"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukrainian critical infrastructure",
      "target_type": "identity"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukraine",
      "target_type": "location"
    },
    {
      "source_name": "Gamaredon",
      "source_type": "threat-actor",
      "relationship_type": "originates-from",
      "target_name": "Russia",
      "target_type": "location"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "NTFS File Attributes",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Windows Management Instrumentation",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Visual Basic",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Automated Collection",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Exfiltration Over C2 Channel",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Worm",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "wscript.exe",
      "target_type": "tool"
    },
    {
      "source_name": "Execution from NTFS Alternate Data Stream",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "VBScript Worm",
      "target_type": "malware"
    },
    {
      "source_name": "Enumeration of NTFS Alternate Data Streams via Dir",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "NTFS File Attributes",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Suspicious Parent for Windows Script Host",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "VBScript Worm",
      "target_type": "malware"
    },
    {
      "source_name": "Ukrainian government organizations",
      "source_type": "identity",
      "relationship_type": "located-at",
      "target_name": "Ukraine",
      "target_type": "location"
    },
    {
      "source_name": "Ukrainian military",
      "source_type": "identity",
      "relationship_type": "located-at",
      "target_name": "Ukraine",
      "target_type": "location"
    },
    {
      "source_name": "Ukrainian critical infrastructure",
      "source_type": "identity",
      "relationship_type": "located-at",
      "target_name": "Ukraine",
      "target_type": "location"
    }
  ]
}