# stix/hola-browser-compromised-in-supply-chain-attack-delivering-cryptominer-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-06-09T15:10:57.686Z",
    "article_slug": "hola-browser-compromised-in-supply-chain-attack-delivering-cryptominer",
    "article_date": "2026-06-09",
    "token_usage": {
      "input_tokens": 143461,
      "output_tokens": 3135,
      "cached_tokens": 137782,
      "total_tokens": 146596,
      "cache_hit": true,
      "estimated_cost_usd": 0.055671
    }
  },
  "report_name": "Hola Browser Compromised in Supply Chain Attack Delivering Cryptominer",
  "threat_actors": [],
  "attack_patterns": [
    {
      "name": "Compromise Software Supply Chain",
      "mitre_id": "T1195.002",
      "description": "Attackers compromised the software distribution pipeline for Hola Browser, injecting a malicious cryptominer into the installation package for version 1.251.91.0."
    },
    {
      "name": "Resource Hijacking",
      "mitre_id": "T1496",
      "description": "A Monero (XMR) cryptominer was installed on user systems to consume CPU cycles and electricity for the attacker's financial gain, a practice known as cryptojacking."
    },
    {
      "name": "Windows Service",
      "mitre_id": "T1543.003",
      "description": "The malware established persistence by creating an auto-starting Windows service named 'hola_monitor_svc'."
    },
    {
      "name": "Disable or Modify Tools",
      "mitre_id": "T1562.001",
      "description": "The malware added an exclusion for its executable (HolaMonitorService.exe) in Windows Defender to evade detection."
    }
  ],
  "malware": [
    {
      "name": "XMRig Cryptominer",
      "mitre_id": "",
      "type": "resource-exploitation",
      "is_family": true,
      "aliases": [
        "me.exe",
        "HolaMonitorService.exe"
      ],
      "description": "A Monero (XMR) cryptominer based on XMRig, distributed via a supply chain attack on Hola Browser. The initial unsigned dropper was named 'me.exe' and it persisted by copying itself to 'HolaMonitorService.exe' and creating a Windows service.",
      "capabilities": [
        "persists-after-system-reboot",
        "degrades-security-software",
        "compromises-system-availability"
      ]
    }
  ],
  "tools": [
    {
      "name": "Hola Browser",
      "mitre_id": "",
      "aliases": [],
      "description": "The Windows version of the Hola Browser (version 1.251.91.0) was compromised and used as the delivery vehicle for a cryptominer in a supply chain attack.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "Hola",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The software vendor whose Windows browser application, Hola Browser, was compromised in a supply chain attack to distribute a cryptominer."
    }
  ],
  "locations": [],
  "infrastructure": [
    {
      "name": "Monero Mining Pool Infrastructure",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "The cryptominer malware connected to Monero (XMR) mining pools to receive work and submit results. Specific pool addresses were not disclosed in the article."
    }
  ],
  "observables": [
    {
      "type": "file",
      "name": "me.exe"
    },
    {
      "type": "file",
      "name": "HolaMonitorService.exe"
    },
    {
      "type": "windows-registry-key",
      "value": "HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\hola_monitor_svc"
    }
  ],
  "indicators": [
    {
      "name": "Hola Browser Cryptominer File Artifacts",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[file:name IN ('me.exe', 'HolaMonitorService.exe')]",
      "x_netsecio_sigma_rule": "title: Hola Browser Cryptominer File Artifacts\nid: 2f7b1d9c-1a8e-4f5c-8d3a-9c4b2e1a0f8d\nstatus: experimental\ndescription: Detects the presence of malicious files associated with the Hola Browser supply chain attack.\nreferences:\n  - https://cyberinsider.com/hola-browser-supply-chain-breach-delivered-crypto-miner-to-users/\nauthor: cyber.netsecops.io\ndate: 2026-06-09\ntags:\n  - attack.t1195.002\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith:\n      - '\\me.exe'\n      - '\\HolaMonitorService.exe'\n  condition: selection",
      "pattern": "",
      "description": "Detects malicious file names associated with the Hola Browser cryptominer payload."
    },
    {
      "name": "Hola Browser Cryptominer Persistence Registry Key",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[windows-registry-key:key = 'HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\hola_monitor_svc']",
      "x_netsecio_sigma_rule": "title: Hola Browser Cryptominer Persistence Registry Key\nid: 8c3a2b1d-9e4f-4a7b-8d3a-9c4b2e1a0f8d\nstatus: experimental\ndescription: Detects the creation of the persistence registry key for the Hola Browser cryptominer service.\nreferences:\n  - https://cyberinsider.com/hola-browser-supply-chain-breach-delivered-crypto-miner-to-users/\nauthor: cyber.netsecops.io\ndate: 2026-06-09\ntags:\n  - attack.t1543.003\nlevel: high\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains: '\\SYSTEM\\CurrentControlSet\\Services\\hola_monitor_svc'\n  condition: selection",
      "pattern": "",
      "description": "Detects the registry key created by the malware for service-based persistence."
    },
    {
      "name": "Hola Browser Cryptominer Service Creation",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Hola Browser Cryptominer Service Creation\nid: 5a9b8c7d-3e2f-4a1b-9c8d-7e6f5d4c3b2a\nstatus: experimental\ndescription: Detects the creation of the malicious 'hola_monitor_svc' Windows service used for persistence by the cryptominer.\nreferences:\n  - https://cyberinsider.com/hola-browser-supply-chain-breach-delivered-crypto-miner-to-users/\nauthor: cyber.netsecops.io\ndate: 2026-06-09\ntags:\n  - attack.t1543.003\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith:\n      - '\\sc.exe'\n      - '\\powershell.exe'\n    CommandLine|contains|all:\n      - 'create'\n      - 'hola_monitor_svc'\n  condition: selection",
      "description": "Detects the creation of the 'hola_monitor_svc' Windows service used for persistence."
    },
    {
      "name": "Hola Browser Cryptominer Windows Defender Exclusion",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Hola Browser Cryptominer Windows Defender Exclusion\nid: 1d2e3f4a-5b6c-4d7e-8f9a-0b1c2d3e4f5a\nstatus: experimental\ndescription: Detects the command used by the Hola Browser cryptominer to add an exclusion for itself in Windows Defender.\nreferences:\n  - https://cyberinsider.com/hola-browser-supply-chain-breach-delivered-crypto-miner-to-users/\nauthor: cyber.netsecops.io\ndate: 2026-06-09\ntags:\n  - attack.t1562.001\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\powershell.exe'\n    CommandLine|contains|all:\n      - 'Add-MpPreference'\n      - '-ExclusionPath'\n      - 'HolaMonitorService.exe'\n  condition: selection",
      "description": "Detects the PowerShell command used to add a Windows Defender exclusion for the malware."
    }
  ],
  "relationships": [
    {
      "source_name": "XMRig Cryptominer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Compromise Software Supply Chain",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "XMRig Cryptominer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Resource Hijacking",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "XMRig Cryptominer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Windows Service",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "XMRig Cryptominer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Disable or Modify Tools",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Hola Browser",
      "source_type": "tool",
      "relationship_type": "delivers",
      "target_name": "XMRig Cryptominer",
      "target_type": "malware"
    },
    {
      "source_name": "XMRig Cryptominer",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "Monero Mining Pool Infrastructure",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Hola Browser Cryptominer File Artifacts",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "XMRig Cryptominer",
      "target_type": "malware"
    },
    {
      "source_name": "Hola Browser Cryptominer File Artifacts",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "me.exe",
      "target_type": "file"
    },
    {
      "source_name": "Hola Browser Cryptominer File Artifacts",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HolaMonitorService.exe",
      "target_type": "file"
    },
    {
      "source_name": "Hola Browser Cryptominer Persistence Registry Key",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "XMRig Cryptominer",
      "target_type": "malware"
    },
    {
      "source_name": "Hola Browser Cryptominer Persistence Registry Key",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\hola_monitor_svc",
      "target_type": "windows-registry-key"
    },
    {
      "source_name": "Hola Browser Cryptominer Service Creation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Windows Service",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Hola Browser Cryptominer Windows Defender Exclusion",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Disable or Modify Tools",
      "target_type": "attack-pattern"
    }
  ]
}