# stix/huntress-uncovers-adware-campaign-exposing-25000-systems-to-supply-chain-attack-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T16:52:31.958Z",
    "article_slug": "huntress-uncovers-adware-campaign-exposing-25000-systems-to-supply-chain-attack",
    "article_date": "2026-04-15",
    "token_usage": {
      "input_tokens": 141624,
      "output_tokens": 4023,
      "cached_tokens": 136033,
      "total_tokens": 145647,
      "cache_hit": true,
      "estimated_cost_usd": 0.064223
    }
  },
  "report_name": "Huntress Uncovers Adware from Dragon Boss Solutions That Disabled AV and Exposed 25,000+ Systems to Trivial Supply Chain Attack",
  "threat_actors": [
    {
      "name": "Dragon Boss Solutions",
      "mitre_id": "",
      "aliases": [],
      "description": "A UAE-based entity that digitally signed adware which disabled security products and created a significant supply chain vulnerability by using an unregistered domain for updates. The software was found on over 25,000 systems.",
      "goals": "Initial goal appears to be typical adware/PUP deployment, but its aggressive capabilities suggest a potential for more malicious follow-on activities.",
      "primary_motivation": "personal-gain",
      "resource_level": "organization",
      "sophistication": "intermediate"
    }
  ],
  "attack_patterns": [
    {
      "name": "Disable or Modify Tools",
      "mitre_id": "T1562.001",
      "description": "A PowerShell script with SYSTEM privileges was used to systematically disable a wide range of cybersecurity products by terminating their processes."
    },
    {
      "name": "Impair Defenses: Disable or Modify System Firewall",
      "mitre_id": "T1562.004",
      "description": "The malware modified the hosts file to block communication with update and telemetry servers of security products."
    },
    {
      "name": "Modify Registry",
      "mitre_id": "T1112",
      "description": "The adware deleted registry keys associated with security products to prevent their reinstallation."
    },
    {
      "name": "Event Triggered Execution: Windows Management Instrumentation Event Subscription",
      "mitre_id": "T1546.003",
      "description": "The malware established persistence through Windows Management Instrumentation (WMI) event subscriptions to ensure it survived reboots."
    },
    {
      "name": "Scheduled Task/Job: Scheduled Task",
      "mitre_id": "T1053.005",
      "description": "The adware created scheduled tasks as a method of persistence."
    },
    {
      "name": "PowerShell",
      "mitre_id": "T1059.001",
      "description": "A PowerShell payload was executed with SYSTEM privileges to perform defense evasion and system modification tasks."
    },
    {
      "name": "Subvert Trust Controls: Code Signing",
      "mitre_id": "T1553.002",
      "description": "The malicious adware was digitally signed by 'Dragon Boss Solutions' to appear as a legitimate program and bypass initial security checks."
    }
  ],
  "malware": [
    {
      "name": "Dragon Boss Solutions Adware",
      "mitre_id": "",
      "type": "adware",
      "is_family": true,
      "aliases": [
        "PUP"
      ],
      "description": "Adware signed by Dragon Boss Solutions that acted as a backdoor, disabling antivirus software, establishing persistence, and creating a supply chain vulnerability by attempting to fetch updates from an unregistered domain.",
      "capabilities": [
        "degrades-security-software",
        "persists-after-system-reboot",
        "communicates-with-c2"
      ]
    }
  ],
  "tools": [
    {
      "name": "PowerShell",
      "mitre_id": "S0194",
      "aliases": [],
      "description": "Used with SYSTEM privileges to disable security products, block their update servers, and delete registry keys.",
      "tool_types": [
        "exploitation"
      ]
    }
  ],
  "vulnerabilities": [
    {
      "name": "Dangling DNS Record for Update Domain",
      "cve": "",
      "description": "The adware attempted to fetch updates from the hardcoded domain 'chromsterabrowser.com', which was unregistered. This created a dangling domain vulnerability, allowing for a trivial supply chain hijack via domain takeover."
    }
  ],
  "campaigns": [],
  "identities": [
    {
      "name": "Universities and colleges",
      "identity_class": "class",
      "sectors": [
        "education"
      ],
      "description": "A class of victims affected by the Dragon Boss Solutions adware, with 221 identified institutions."
    },
    {
      "name": "Operational Technology (OT) networks",
      "identity_class": "class",
      "sectors": [
        "infrastructure",
        "utilities"
      ],
      "description": "A class of victims affected by the Dragon Boss Solutions adware, with 41 identified OT networks, including electric utilities."
    },
    {
      "name": "Government entities",
      "identity_class": "class",
      "sectors": [
        "government"
      ],
      "description": "A class of victims affected by the Dragon Boss Solutions adware, with 35 identified government entities."
    },
    {
      "name": "Healthcare organizations",
      "identity_class": "class",
      "sectors": [
        "healthcare"
      ],
      "description": "A class of victims affected by the Dragon Boss Solutions adware, with 3 identified healthcare organizations."
    }
  ],
  "locations": [
    {
      "name": "United Arab Emirates",
      "region": "western-asia"
    }
  ],
  "infrastructure": [
    {
      "name": "Dragon Boss Adware Update Server",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "An update server infrastructure that the adware attempted to contact via the unregistered domain 'chromsterabrowser.com'. This created a major supply chain risk."
    },
    {
      "name": "Dragon Boss Adware C2 Server",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "A command and control server associated with the campaign, using the domain 'worldwidewebframework3.com'."
    }
  ],
  "observables": [
    {
      "type": "domain-name",
      "value": "chromsterabrowser.com"
    },
    {
      "type": "domain-name",
      "value": "worldwidewebframework3.com"
    }
  ],
  "indicators": [
    {
      "name": "Dragon Boss Solutions Adware C2/Update Domains",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[domain-name:value IN ('chromsterabrowser.com', 'worldwidewebframework3.com')]",
      "x_netsecio_sigma_rule": "title: Dragon Boss Solutions Adware C2 DNS Query\nid: 7f5b1a3c-8d2e-4f9a-9c1b-0a1b2c3d4e5f\nstatus: experimental\ndescription: Detects DNS queries to domains associated with the Dragon Boss Solutions adware campaign.\nreferences:\n  - https://www.huntress.com/blog/pups-grow-fangs-dragon-boss-solutions-supply-chain-risk\nauthor: cyber.netsecops.io\ndate: 2026-04-15\ntags:\n  - attack.t1071.004\n  - attack.command_and_control\nlevel: high\nlogsource:\n  category: dns_query\n  product: windows\ndetection:\n  selection:\n    QueryName|contains:\n      - 'chromsterabrowser.com'\n      - 'worldwidewebframework3.com'\n  condition: selection",
      "pattern": "",
      "description": "Detects network connections to the C2 and update domains used by the Dragon Boss Solutions adware."
    },
    {
      "name": "Dragon Boss Adware Disables Windows Defender via PowerShell",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Dragon Boss Adware Disables Windows Defender via PowerShell\nid: 8a9b0c1d-2e3f-4a5b-6c7d-8e9f0a1b2c3d\nstatus: experimental\ndescription: Detects the execution of a PowerShell command used by Dragon Boss Solutions adware to disable Windows Defender real-time monitoring.\nreferences:\n  - https://www.huntress.com/blog/pups-grow-fangs-dragon-boss-solutions-supply-chain-risk\nauthor: cyber.netsecops.io\ndate: 2026-04-15\ntags:\n  - attack.t1562.001\n  - attack.t1059.001\n  - attack.defense_evasion\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\powershell.exe'\n    CommandLine|contains|all:\n      - 'Set-MpPreference'\n      - '-DisableRealtimeMonitoring'\n      - '$true'\n  condition: selection",
      "description": "Detects the specific PowerShell command used by the adware to disable Windows Defender's real-time monitoring."
    },
    {
      "name": "Dragon Boss Adware Hosts File Modification",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Dragon Boss Adware Hosts File Modification\nid: 9b0c1d2e-3f4a-5b6c-7d8e-9f0a1b2c3d4e\nstatus: experimental\ndescription: Detects modification of the hosts file, a technique used by Dragon Boss Solutions adware to block security product updates.\nreferences:\n  - https://www.huntress.com/blog/pups-grow-fangs-dragon-boss-solutions-supply-chain-risk\nauthor: cyber.netsecops.io\ndate: 2026-04-15\ntags:\n  - attack.t1562.004\n  - attack.defense_evasion\nlevel: medium\nlogsource:\n  category: file_event\n  product: windows\ndetection:\n  selection:\n    TargetFilename|endswith: '\\drivers\\etc\\hosts'\n  filter_legit:\n    Image: \n      - 'C:\\Windows\\System32\\notepad.exe'\n  condition: selection and not filter_legit",
      "description": "Detects modifications to the system's hosts file, a technique used by the adware to block security vendor communications."
    }
  ],
  "relationships": [
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Dragon Boss Solutions Adware",
      "target_type": "malware"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "PowerShell",
      "target_type": "tool"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "originates-from",
      "target_name": "United Arab Emirates",
      "target_type": "location"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Universities and colleges",
      "target_type": "identity"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Operational Technology (OT) networks",
      "target_type": "identity"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Government entities",
      "target_type": "identity"
    },
    {
      "source_name": "Dragon Boss Solutions",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Healthcare organizations",
      "target_type": "identity"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Disable or Modify Tools",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Impair Defenses: Disable or Modify System Firewall",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Modify Registry",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Event Triggered Execution: Windows Management Instrumentation Event Subscription",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Scheduled Task/Job: Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "PowerShell",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Subvert Trust Controls: Code Signing",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "exploits",
      "target_name": "Dangling DNS Record for Update Domain",
      "target_type": "vulnerability"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "communicates-with",
      "target_name": "Dragon Boss Adware Update Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Dragon Boss Solutions Adware",
      "source_type": "malware",
      "relationship_type": "communicates-with",
      "target_name": "Dragon Boss Adware C2 Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Dragon Boss Solutions Adware C2/Update Domains",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Dragon Boss Solutions Adware",
      "target_type": "malware"
    },
    {
      "source_name": "Dragon Boss Solutions Adware C2/Update Domains",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "chromsterabrowser.com",
      "target_type": "domain-name"
    },
    {
      "source_name": "Dragon Boss Solutions Adware C2/Update Domains",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "worldwidewebframework3.com",
      "target_type": "domain-name"
    },
    {
      "source_name": "Dragon Boss Adware Disables Windows Defender via PowerShell",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Disable or Modify Tools",
      "target_type": "attack-pattern"
    }
  ]
}