# stix/jscrambler-npm-package-compromised-supply-chain-attack-infostealer-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-07-12T13:35:16.909Z",
    "article_slug": "jscrambler-npm-package-compromised-supply-chain-attack-infostealer",
    "article_date": "2026-07-12",
    "token_usage": {
      "input_tokens": 146626,
      "output_tokens": 4311,
      "cached_tokens": 140082,
      "total_tokens": 150937,
      "cache_hit": true,
      "estimated_cost_usd": 0.0688
    }
  },
  "report_name": "Jscrambler NPM Package Hijacked in Supply Chain Attack, Deploys Infostealer Targeting Cloud & AI Dev Secrets",
  "threat_actors": [],
  "attack_patterns": [
    {
      "name": "Compromise Software Supply Chain",
      "mitre_id": "T1195.002",
      "description": "The attacker compromised a maintainer's npm publishing token to publish malicious versions of the 'jscrambler' package to the public npm registry."
    },
    {
      "name": "Command and Scripting Interpreter",
      "mitre_id": "T1059",
      "description": "The malicious npm package used a 'preinstall' hook in the package.json file to automatically execute a setup script that deployed an infostealer payload."
    },
    {
      "name": "Cloud Instance Metadata API",
      "mitre_id": "T1552.005",
      "description": "The infostealer was designed to search for and exfiltrate cloud credentials for AWS, Azure, and Google Cloud from developer workstations and CI/CD environments."
    },
    {
      "name": "Credentials In Files",
      "mitre_id": "T1552.001",
      "description": "The infostealer searched for and exfiltrated cryptocurrency wallet data from MetaMask, Phantom, and Exodus, as well as configuration files for AI coding tools."
    },
    {
      "name": "Scheduled Task",
      "mitre_id": "T1053.005",
      "description": "On Windows systems, the infostealer attempted to establish persistence by creating a hidden scheduled task."
    },
    {
      "name": "Launch Agent",
      "mitre_id": "T1543.001",
      "description": "On macOS systems, the infostealer used a LaunchAgent to establish persistence."
    }
  ],
  "malware": [
    {
      "name": "Rust Infostealer",
      "mitre_id": "",
      "type": "infostealer",
      "is_family": false,
      "aliases": [],
      "description": "A cross-platform infostealer written in Rust, distributed via malicious 'jscrambler' npm packages. It was designed to harvest sensitive developer credentials from cloud providers (AWS, Azure, GCP), cryptocurrency wallets, and AI coding tools.",
      "capabilities": [
        "steals-authentication-credentials",
        "persists-after-system-reboot"
      ]
    }
  ],
  "tools": [
    {
      "name": "npm",
      "mitre_id": "",
      "aliases": [],
      "description": "The npm package registry was used as the distribution vector for the malicious 'jscrambler' package versions containing the infostealer."
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "Jscrambler",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The vendor whose 'jscrambler' npm package was hijacked in a software supply chain attack."
    },
    {
      "name": "npm",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The public package registry that hosted and distributed the malicious versions of the 'jscrambler' package."
    },
    {
      "name": "AWS",
      "identity_class": "organization",
      "sectors": [
        "technology",
        "infrastructure"
      ],
      "description": "Cloud provider whose developer credentials were targeted by the Rust infostealer."
    },
    {
      "name": "Azure",
      "identity_class": "organization",
      "sectors": [
        "technology",
        "infrastructure"
      ],
      "description": "Cloud provider whose developer credentials were targeted by the Rust infostealer."
    },
    {
      "name": "Google Cloud",
      "identity_class": "organization",
      "sectors": [
        "technology",
        "infrastructure"
      ],
      "description": "Cloud provider whose developer credentials were targeted by the Rust infostealer."
    },
    {
      "name": "MetaMask",
      "identity_class": "organization",
      "sectors": [
        "financial-services",
        "technology"
      ],
      "description": "Cryptocurrency wallet whose data was targeted by the Rust infostealer."
    },
    {
      "name": "Phantom",
      "identity_class": "organization",
      "sectors": [
        "financial-services",
        "technology"
      ],
      "description": "Cryptocurrency wallet whose data was targeted by the Rust infostealer."
    },
    {
      "name": "Exodus",
      "identity_class": "organization",
      "sectors": [
        "financial-services",
        "technology"
      ],
      "description": "Cryptocurrency wallet whose data was targeted by the Rust infostealer."
    },
    {
      "name": "Claude Desktop",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "AI coding tool whose configuration files and credentials were targeted by the Rust infostealer."
    },
    {
      "name": "Cursor",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "AI coding tool whose configuration files and credentials were targeted by the Rust infostealer."
    },
    {
      "name": "VS Code",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "Development environment whose configuration files were targeted by the Rust infostealer."
    }
  ],
  "locations": [],
  "infrastructure": [],
  "observables": [
    {
      "type": "file",
      "name": "intro.js"
    }
  ],
  "indicators": [
    {
      "name": "Jscrambler Supply Chain Attack Payload",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[file:name = 'intro.js' AND file:parent_directory_ref.path MATCHES '.*(\\\\|/)node_modules(\\\\|/)jscrambler$']",
      "x_netsecio_sigma_rule": "title: Jscrambler Supply Chain Attack Payload\nid: 4f5a6b7c-8d9e-4f1a-b2c3-d4e5f6a7b8c9\nstatus: experimental\ndescription: Detects the creation of the malicious 'intro.js' payload from the compromised jscrambler npm package.\nauthor: cyber.netsecops.io\ndate: 2026-07-12\nreferences:\n  - https://socket.dev/blog/jscrambler-supply-chain-attack\ntags:\n  - attack.t1195.002\nlevel: critical\nlogsource:\n  category: file_event\n  product: windows\ndetection:\n  selection:\n    TargetFilename|endswith: '\\node_modules\\jscrambler\\intro.js'\n  condition: selection",
      "pattern": "",
      "description": "Detects the presence of the malicious payload file 'intro.js' dropped by compromised versions of the 'jscrambler' npm package."
    },
    {
      "name": "Suspicious Jscrambler NPM Package Installation",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Jscrambler NPM Package Installation\nid: 1a2b3c4d-5e6f-4a7b-8c9d-0e1f2a3b4c5d\nstatus: experimental\ndescription: Detects the installation of specific malicious versions of the 'jscrambler' npm package via command line arguments.\nauthor: cyber.netsecops.io\ndate: 2026-07-12\nreferences:\n  - https://socket.dev/blog/jscrambler-supply-chain-attack\ntags:\n  - attack.t1195.002\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_base:\n    Image|endswith:\n      - '\\npm.cmd'\n      - '\\node.exe'\n    CommandLine|contains|all:\n      - 'install'\n      - 'jscrambler@'\n  selection_versions:\n    CommandLine|contains:\n      - '8.14.0'\n      - '8.16.0'\n      - '8.17.0'\n      - '8.18.0'\n      - '8.20.0'\n  condition: all of selection_*\n",
      "description": "Detects command-line execution of 'npm install' for known malicious versions of the 'jscrambler' package."
    },
    {
      "name": "Suspicious Scheduled Task Creation by Node Process",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Scheduled Task Creation by Node Process\nid: 6f7a8b9c-0d1e-4f2a-b3c4-d5e6f7a8b9c0\nstatus: experimental\ndescription: Detects a scheduled task being created by a node.js or npm process, which could indicate persistence from a malicious npm package.\nauthor: cyber.netsecops.io\ndate: 2026-07-12\nreferences:\n  - https://socket.dev/blog/jscrambler-supply-chain-attack\ntags:\n  - attack.t1053.005\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    ParentImage|endswith:\n      - '\\node.exe'\n      - '\\npm.cmd'\n    Image|endswith: '\\schtasks.exe'\n    CommandLine|contains: '/create'\n  condition: selection\n",
      "description": "Detects the creation of a Windows Scheduled Task by a node or npm process, a TTP used by the malicious jscrambler package for persistence."
    },
    {
      "name": "Suspicious LaunchAgent Creation by Node Process",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious LaunchAgent Creation by Node Process\nid: 7a8b9c0d-1e2f-4a3b-b4c5-d6e7f8a9b0c1\nstatus: experimental\ndescription: Detects a LaunchAgent .plist file being created by a node.js process, which could indicate persistence from a malicious npm package on macOS.\nauthor: cyber.netsecops.io\ndate: 2026-07-12\nreferences:\n  - https://socket.dev/blog/jscrambler-supply-chain-attack\ntags:\n  - attack.t1543.001\nlevel: medium\nlogsource:\n  category: file_event\n  product: macos\ndetection:\n  selection:\n    TargetFilename|contains: '/Library/LaunchAgents/'\n    TargetFilename|endswith: '.plist'\n    Image|endswith: '/node'\n  condition: selection\n",
      "description": "Detects the creation of a macOS LaunchAgent by a node process, a TTP used by the malicious jscrambler package for persistence."
    }
  ],
  "relationships": [
    {
      "source_name": "Compromise Software Supply Chain",
      "source_type": "attack-pattern",
      "relationship_type": "delivers",
      "target_name": "Rust Infostealer",
      "target_type": "malware"
    },
    {
      "source_name": "Compromise Software Supply Chain",
      "source_type": "attack-pattern",
      "relationship_type": "targets",
      "target_name": "Jscrambler",
      "target_type": "identity"
    },
    {
      "source_name": "Compromise Software Supply Chain",
      "source_type": "attack-pattern",
      "relationship_type": "uses",
      "target_name": "npm",
      "target_type": "tool"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Command and Scripting Interpreter",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Cloud Instance Metadata API",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Credentials In Files",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Launch Agent",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "AWS",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Azure",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Google Cloud",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "MetaMask",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Phantom",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Exodus",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Claude Desktop",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "Cursor",
      "target_type": "identity"
    },
    {
      "source_name": "Rust Infostealer",
      "source_type": "malware",
      "relationship_type": "targets",
      "target_name": "VS Code",
      "target_type": "identity"
    },
    {
      "source_name": "Jscrambler Supply Chain Attack Payload",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "intro.js",
      "target_type": "file"
    },
    {
      "source_name": "Jscrambler Supply Chain Attack Payload",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Rust Infostealer",
      "target_type": "malware"
    },
    {
      "source_name": "Suspicious Jscrambler NPM Package Installation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Compromise Software Supply Chain",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Suspicious Scheduled Task Creation by Node Process",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Suspicious LaunchAgent Creation by Node Process",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Launch Agent",
      "target_type": "attack-pattern"
    }
  ]
}