# stix/megalodon-supply-chain-attack-compromises-over-5500-github-repositories-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-05-30T18:07:11.052Z",
    "article_slug": "megalodon-supply-chain-attack-compromises-over-5500-github-repositories",
    "article_date": "2026-05-30",
    "token_usage": {
      "input_tokens": 144995,
      "output_tokens": 3968,
      "cached_tokens": 137782,
      "total_tokens": 148963,
      "cache_hit": true,
      "estimated_cost_usd": 0.065919
    }
  },
  "report_name": "Massive 'Megalodon' Supply Chain Attack Compromises 5,500+ GitHub Repos to Steal Cloud Credentials",
  "threat_actors": [
    {
      "name": "TeamPCP",
      "mitre_id": "",
      "aliases": [],
      "description": "Threat actor attributed to the 'Megalodon' and 'Mini Shai-Hulud' supply chain attack campaigns, targeting GitHub repositories and package ecosystems (npm, PyPI) to steal cloud credentials and other secrets from CI/CD environments.",
      "goals": "Mass theft of sensitive credentials and secrets (AWS, GCP, Azure, API tokens, SSH keys) from CI/CD environments.",
      "sophistication": "advanced"
    }
  ],
  "attack_patterns": [
    {
      "name": "Compromise Software Dependencies and Development Tools",
      "mitre_id": "T1195.001",
      "description": "The core of the attack was compromising GitHub repositories and injecting malicious GitHub Actions, a key development tool in the software supply chain."
    },
    {
      "name": "Default Accounts",
      "mitre_id": "T1078.001",
      "description": "Attackers likely gained initial access to the GitHub repositories by using compromised developer credentials or tokens, which may have included default credentials."
    },
    {
      "name": "Unsecured Credentials",
      "mitre_id": "T1552",
      "description": "The primary objective of the Megalodon attack was to steal cloud credentials (AWS, GCP, Azure), API tokens, and SSH keys stored as secrets or environment variables within the GitHub Actions CI/CD environment."
    },
    {
      "name": "Automated Collection",
      "mitre_id": "T1119",
      "description": "The malicious GitHub Action workflow automatically collected secrets and environment variables from the CI/CD environment for exfiltration."
    },
    {
      "name": "Scheduled Task/Job",
      "mitre_id": "T1053",
      "description": "The malicious GitHub Action workflow was configured to run automatically on 'push' and 'pull_request' events, acting as a scheduled job to execute the attacker's code within the trusted context of the repository's runner."
    },
    {
      "name": "Exfiltration Over C2 Channel",
      "mitre_id": "T1041",
      "description": "The malicious workflow exfiltrated the stolen credentials and secrets using tools like curl or wget to an external attacker-controlled domain."
    }
  ],
  "malware": [
    {
      "name": "Malicious GitHub Action Workflow",
      "mitre_id": "",
      "type": "infostealer",
      "is_family": true,
      "aliases": [],
      "description": "Malicious GitHub Action workflows injected into repositories by the Megalodon campaign. These workflows were configured to trigger on 'push' and 'pull_request' events to exfiltrate secrets, environment variables, and install dormant backdoors.",
      "capabilities": [
        "steals-authentication-credentials",
        "exfiltrates-data",
        "communicates-with-c2"
      ]
    }
  ],
  "tools": [
    {
      "name": "base64",
      "mitre_id": "",
      "aliases": [],
      "description": "Used within the malicious GitHub Action workflow to obfuscate exfiltrated environment variables and secrets before exfiltration."
    },
    {
      "name": "curl",
      "mitre_id": "",
      "aliases": [],
      "description": "Used within the malicious GitHub Action workflow to exfiltrate stolen credentials and secrets to an attacker-controlled domain."
    },
    {
      "name": "wget",
      "mitre_id": "",
      "aliases": [],
      "description": "Mentioned as a potential tool used alongside curl within the malicious GitHub Action workflow to exfiltrate stolen credentials and secrets."
    }
  ],
  "vulnerabilities": [],
  "campaigns": [
    {
      "name": "Megalodon",
      "description": "A large-scale supply chain attack campaign that compromised over 5,500 GitHub repositories by injecting malicious GitHub Action workflows to exfiltrate cloud credentials and other secrets. This was the second wave of a broader campaign attributed to TeamPCP.",
      "first_seen": "2026-05-18T00:00:00Z",
      "objective": "Exfiltrate secrets from CI/CD environments of over 5,500 GitHub repositories."
    },
    {
      "name": "Mini Shai-Hulud",
      "description": "The first wave of a supply chain attack campaign, preceding 'Megalodon'. It was a self-propagating worm that compromised 172 packages on npm and PyPI.",
      "first_seen": "2026-04-29T00:00:00Z",
      "objective": "Compromise npm and PyPI packages."
    }
  ],
  "identities": [
    {
      "name": "GitHub",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The code hosting platform where over 5,500 repositories were compromised by the Megalodon supply chain attack."
    },
    {
      "name": "npm",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "Package registry for JavaScript, targeted by the 'Mini Shai-Hulud' worm, the first wave of the campaign."
    },
    {
      "name": "PyPI",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The Python Package Index, targeted by the 'Mini Shai-Hulud' worm, the first wave of the campaign."
    },
    {
      "name": "Amazon Web Services",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "Cloud provider whose credentials were a primary target for exfiltration in the Megalodon attack."
    },
    {
      "name": "Google Cloud Platform",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "Cloud provider whose access tokens were a primary target for exfiltration in the Megalodon attack."
    },
    {
      "name": "Microsoft Azure",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "Cloud provider whose credentials were a primary target for exfiltration in the Megalodon attack."
    }
  ],
  "locations": [],
  "infrastructure": [
    {
      "name": "Megalodon Exfiltration Endpoint",
      "infrastructure_types": [
        "command-and-control",
        "exfiltration"
      ],
      "description": "An external, non-standard domain used as an endpoint for exfiltrating stolen credentials via curl or wget from compromised GitHub Action workflows."
    }
  ],
  "observables": [],
  "indicators": [
    {
      "name": "Suspicious GitHub Workflow Modification",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious GitHub Workflow Modification\nid: 2b8c1f9a-3d7e-4b0f-8c6a-9e2d1f8a7b3c\nstatus: experimental\ndescription: Monitors for unexpected or unauthorized commits that add or modify files in the .github/workflows/ directory.\nauthor: cyber.netsecops.io\ndate: 2026-05-30\ntags:\n  - attack.t1195.001\nlevel: medium\nlogsource:\n  category: file_event\ndetection:\n  selection:\n    TargetFilename|contains: '/.github/workflows/'\n  condition: selection",
      "description": "Detects modifications to GitHub workflow files, which was the primary vector for the Megalodon attack."
    },
    {
      "name": "GitHub Workflow Environment Variable Exfiltration",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: GitHub Workflow Environment Variable Exfiltration\nid: 8e4d2c1b-9a0f-4d7e-8c6a-3d7e1f8a7b3c\nstatus: experimental\ndescription: Detects GitHub Action workflow steps that pipe environment variables to encoding commands like base64, a common technique to obfuscate exfiltrated data.\nauthor: cyber.netsecops.io\ndate: 2026-05-30\ntags:\n  - attack.t1119\n  - attack.t1041\nlevel: high\nlogsource:\n  category: process_creation\ndetection:\n  selection:\n    CommandLine|contains: 'env | base64'\n  condition: selection",
      "description": "Detects a command pattern used in the Megalodon attack to collect and obfuscate environment variables for exfiltration."
    },
    {
      "name": "GitHub Workflow Data Exfiltration via cURL",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: GitHub Workflow Data Exfiltration via cURL\nid: f1a2b3c4-d5e6-4a7b-8c9d-0e1f2a3b4c5d\nstatus: experimental\ndescription: Detects GitHub Action workflow run commands that use curl or wget to send data to an external, non-standard domain.\nauthor: cyber.netsecops.io\ndate: 2026-05-30\ntags:\n  - attack.t1041\nlevel: high\nlogsource:\n  category: process_creation\ndetection:\n  selection_curl:\n    Image|endswith:\n      - '/curl'\n      - '/wget'\n    CommandLine|contains|all:\n      - '-X POST'\n      - '-d @-'\n  condition: selection_curl",
      "description": "Detects the use of curl or wget with specific parameters to exfiltrate data, a pattern observed in the Megalodon attack."
    }
  ],
  "relationships": [
    {
      "source_name": "TeamPCP",
      "source_type": "threat-actor",
      "relationship_type": "attributed-to",
      "target_name": "Megalodon",
      "target_type": "campaign"
    },
    {
      "source_name": "TeamPCP",
      "source_type": "threat-actor",
      "relationship_type": "attributed-to",
      "target_name": "Mini Shai-Hulud",
      "target_type": "campaign"
    },
    {
      "source_name": "TeamPCP",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Malicious GitHub Action Workflow",
      "target_type": "malware"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Compromise Software Dependencies and Development Tools",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Unsecured Credentials",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Default Accounts",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Scheduled Task/Job",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "GitHub",
      "target_type": "identity"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "Amazon Web Services",
      "target_type": "identity"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "Google Cloud Platform",
      "target_type": "identity"
    },
    {
      "source_name": "Megalodon",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "Microsoft Azure",
      "target_type": "identity"
    },
    {
      "source_name": "Mini Shai-Hulud",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "npm",
      "target_type": "identity"
    },
    {
      "source_name": "Mini Shai-Hulud",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "PyPI",
      "target_type": "identity"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Automated Collection",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Exfiltration Over C2 Channel",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "base64",
      "target_type": "tool"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "curl",
      "target_type": "tool"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "wget",
      "target_type": "tool"
    },
    {
      "source_name": "Malicious GitHub Action Workflow",
      "source_type": "malware",
      "relationship_type": "exfiltrates-to",
      "target_name": "Megalodon Exfiltration Endpoint",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Suspicious GitHub Workflow Modification",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Compromise Software Dependencies and Development Tools",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "GitHub Workflow Environment Variable Exfiltration",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Automated Collection",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "GitHub Workflow Data Exfiltration via cURL",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Exfiltration Over C2 Channel",
      "target_type": "attack-pattern"
    }
  ]
}