# stix/new-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-06-25T14:58:11.455Z",
    "article_slug": "new-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke",
    "article_date": "2026-06-25",
    "token_usage": {
      "input_tokens": 146518,
      "output_tokens": 5145,
      "cached_tokens": 140082,
      "total_tokens": 151663,
      "cache_hit": true,
      "estimated_cost_usd": 0.077005
    }
  },
  "report_name": "Stealthy 'Mistic' Backdoor Linked to Ransomware IAB KongTuke in Pre-Attack Campaigns",
  "threat_actors": [
    {
      "name": "KongTuke",
      "mitre_id": "",
      "aliases": [
        "Woodgnat"
      ],
      "description": "A prolific initial access broker (IAB) that compromises corporate networks, establishes persistent access using custom tools like Mistic, and then sells that access to ransomware groups. The group has been active since at least April 2026.",
      "goals": "To gain and maintain persistent, stealthy access to corporate networks for the purpose of selling that access to ransomware operators.",
      "primary_motivation": "personal-gain",
      "sophistication": "advanced"
    },
    {
      "name": "Qilin",
      "mitre_id": "G1003",
      "aliases": [],
      "description": "A major ransomware group that purchases network access from initial access brokers like KongTuke to deploy their ransomware.",
      "primary_motivation": "personal-gain",
      "sophistication": "advanced"
    },
    {
      "name": "Akira",
      "mitre_id": "G1024",
      "aliases": [],
      "description": "A major ransomware group that purchases network access from initial access brokers like KongTuke to deploy their ransomware.",
      "primary_motivation": "personal-gain",
      "sophistication": "advanced"
    },
    {
      "name": "Black Basta",
      "mitre_id": "",
      "aliases": [],
      "description": "A major ransomware group that purchases network access from initial access brokers like KongTuke to deploy their ransomware.",
      "primary_motivation": "personal-gain",
      "sophistication": "advanced"
    }
  ],
  "attack_patterns": [
    {
      "name": "DLL Side-Loading",
      "mitre_id": "T1574.002",
      "description": "The Mistic backdoor was loaded by a legitimate, signed executable (MpExtMs.exe) which loaded a malicious DLL (EndpointDlp.dll) placed in the same directory. This technique is used to bypass application allowlisting and reputation-based security controls."
    },
    {
      "name": "Match Legitimate Resource Name or Location",
      "mitre_id": "T1036.005",
      "description": "The malicious DLL used for side-loading was named 'EndpointDlp.dll' to masquerade as a component of Microsoft's endpoint security suite."
    },
    {
      "name": "File Deletion",
      "mitre_id": "T1070.004",
      "description": "The Mistic backdoor includes a kill switch that allows it to delete itself from the compromised system to hinder investigation."
    },
    {
      "name": "Web Protocols",
      "mitre_id": "T1071.001",
      "description": "The Mistic backdoor communicates with its C2 server over standard web protocols to receive commands and exfiltrate data, blending in with normal network traffic."
    },
    {
      "name": "Registry Run Keys / Startup Folder",
      "mitre_id": "T1547.001",
      "description": "KongTuke establishes persistence for the abused executable (MpExtMs.exe) by creating entries in the 'HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run' registry key."
    },
    {
      "name": "Regsvr32",
      "mitre_id": "T1218.010",
      "description": "KongTuke has been observed using the command 'regsvr32 /s /n /i:<url> scrobj.dll' to execute remote scripts, a technique associated with Squiblydoo."
    }
  ],
  "malware": [
    {
      "name": "Mistic",
      "mitre_id": "",
      "type": "backdoor",
      "is_family": true,
      "aliases": [
        "MLTBackdoor"
      ],
      "description": "A sophisticated and stealthy Windows backdoor used by the KongTuke IAB. It is designed for long-term, low-visibility persistence, operating primarily in memory and using DLL side-loading for execution. It serves as a foothold for future ransomware deployment.",
      "capabilities": [
        "communicates-with-c2",
        "cleans-traces-of-infection",
        "persists-after-system-reboot"
      ]
    },
    {
      "name": "ModeloRAT",
      "mitre_id": "",
      "type": "remote-access-trojan",
      "is_family": true,
      "aliases": [],
      "description": "A remote access trojan deployed by KongTuke, observed in at least one incident alongside the Mistic backdoor."
    },
    {
      "name": "Qilin",
      "mitre_id": "S1242",
      "type": "ransomware",
      "is_family": true,
      "aliases": [],
      "description": "Ransomware deployed by affiliates who purchase access from initial access brokers like KongTuke."
    },
    {
      "name": "Akira",
      "mitre_id": "S1129",
      "type": "ransomware",
      "is_family": true,
      "aliases": [],
      "description": "Ransomware deployed by affiliates who purchase access from initial access brokers like KongTuke."
    },
    {
      "name": "Black Basta",
      "mitre_id": "S1068",
      "type": "ransomware",
      "is_family": true,
      "aliases": [],
      "description": "Ransomware deployed by affiliates who purchase access from initial access brokers like KongTuke."
    }
  ],
  "tools": [
    {
      "name": "MpExtMs.exe",
      "mitre_id": "",
      "aliases": [],
      "description": "A legitimate, signed Microsoft binary that was abused by KongTuke to perform DLL side-loading of the Mistic backdoor.",
      "tool_types": []
    },
    {
      "name": "Regsvr32",
      "mitre_id": "S0108",
      "aliases": [],
      "description": "The threat actor KongTuke uses regsvr32.exe to execute remote scripts as part of their attack chain, a technique known as Squiblydoo.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [
    {
      "name": "Insurance Sector Organizations",
      "identity_class": "class",
      "sectors": [
        "insurance"
      ],
      "description": "Organizations in the insurance sector were targeted by the KongTuke group using the Mistic backdoor."
    },
    {
      "name": "Education Sector Organizations",
      "identity_class": "class",
      "sectors": [
        "education"
      ],
      "description": "Organizations in the education sector were targeted by the KongTuke group using the Mistic backdoor."
    },
    {
      "name": "IT Sector Organizations",
      "identity_class": "class",
      "sectors": [
        "technology"
      ],
      "description": "Organizations in the IT sector were targeted by the KongTuke group using the Mistic backdoor."
    }
  ],
  "locations": [],
  "infrastructure": [
    {
      "name": "Mistic C2 Infrastructure",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "Command and control infrastructure used by the Mistic backdoor to receive commands and exfiltrate data."
    }
  ],
  "observables": [
    {
      "type": "file",
      "name": "EndpointDlp.dll"
    },
    {
      "type": "file",
      "name": "MpExtMs.exe"
    },
    {
      "type": "windows-registry-key",
      "value": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
    }
  ],
  "indicators": [
    {
      "name": "Mistic Backdoor DLL Side-Loading",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Mistic Backdoor DLL Side-Loading Execution\nid: 5e8a1b3c-7d4f-4e9a-8b2c-1d6f0a8b4c5d\nstatus: experimental\ndescription: Detects the execution of a legitimate Microsoft executable (MpExtMs.exe) from a non-standard directory, which is a known TTP for the KongTuke IAB to side-load the Mistic backdoor (EndpointDlp.dll).\nreferences:\n  - https://www.bleepingcomputer.com/news/security/stealthy-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke/\nauthor: cyber.netsecops.io\ndate: 2026-06-25\ntags:\n  - attack.t1574.002\n  - attack.defense_evasion\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\MpExtMs.exe'\n  filter:\n    Image|startswith: 'C:\\Windows\\System32\\'\n  condition: selection and not filter",
      "description": "Detects the execution of MpExtMs.exe from a non-standard directory, a key indicator of the Mistic backdoor's DLL side-loading technique."
    },
    {
      "name": "KongTuke Squiblydoo Execution",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: KongTuke Associated Squiblydoo Command-Line\nid: a2b3c4d5-e6f7-4a8b-9c0d-1e2f3a4b5c6d\nstatus: experimental\ndescription: Detects the use of regsvr32 to execute remote scripts, a technique associated with KongTuke.\nreferences:\n  - https://www.bleepingcomputer.com/news/security/stealthy-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke/\nauthor: cyber.netsecops.io\ndate: 2026-06-25\ntags:\n  - attack.t1218.010\n  - attack.defense_evasion\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\regsvr32.exe'\n    CommandLine|contains|all:\n      - '/s'\n      - '/n'\n      - '/i:'\n      - 'scrobj.dll'\n  condition: selection",
      "description": "Detects the specific command-line pattern used by KongTuke to execute remote scripts via regsvr32, a technique known as Squiblydoo."
    },
    {
      "name": "Mistic Backdoor Persistence via Run Key",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[windows-registry-key:key = 'HKEY_CURRENT_USER\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run'] AND [file:name = 'MpExtMs.exe']",
      "x_netsecio_sigma_rule": "title: Mistic Backdoor Persistence via Run Key\nid: f1e2d3c4-b5a6-4987-8654-321fedcba987\nstatus: experimental\ndescription: Detects the creation of a registry run key pointing to MpExtMs.exe, which is used by KongTuke for Mistic backdoor persistence.\nreferences:\n  - https://www.bleepingcomputer.com/news/security/stealthy-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke/\nauthor: cyber.netsecops.io\ndate: 2026-06-25\ntags:\n  - attack.t1547.001\n  - attack.persistence\nlevel: high\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    EventType: 'SetValue'\n    TargetObject|contains: '\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n    Details|contains: 'MpExtMs.exe'\n  condition: selection",
      "pattern": "",
      "description": "Detects the creation of a persistence entry in the HKCU Run key pointing to the abused MpExtMs.exe executable."
    },
    {
      "name": "Mistic Backdoor Associated File Names",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[file:name IN ('EndpointDlp.dll', 'MpExtMs.exe')]",
      "x_netsecio_sigma_rule": "title: Mistic Backdoor Associated File Creation\nid: 98765432-1fed-cba9-8765-43210fedcba9\nstatus: experimental\ndescription: Detects the presence of files associated with the Mistic backdoor side-loading routine.\nreferences:\n  - https://www.bleepingcomputer.com/news/security/stealthy-mistic-backdoor-linked-to-ransomware-access-broker-kongtuke/\nauthor: cyber.netsecops.io\ndate: 2026-06-25\ntags:\n  - attack.t1574.002\nlevel: medium\nlogsource:\n  category: file_event\n  product: windows\ndetection:\n  selection:\n    TargetFilename|endswith:\n      - '\\EndpointDlp.dll'\n      - '\\MpExtMs.exe'\n  filter:\n    TargetFilename|startswith: 'C:\\Windows\\System32\\'\n  condition: selection and not filter",
      "pattern": "",
      "description": "Detects file names associated with the Mistic backdoor, specifically the malicious DLL and the abused executable, especially when found outside their legitimate directories."
    }
  ],
  "relationships": [
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Mistic",
      "target_type": "malware"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "ModeloRAT",
      "target_type": "malware"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "DLL Side-Loading",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Regsvr32",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "related-to",
      "target_name": "Qilin",
      "target_type": "threat-actor"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "related-to",
      "target_name": "Akira",
      "target_type": "threat-actor"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "related-to",
      "target_name": "Black Basta",
      "target_type": "threat-actor"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Insurance Sector Organizations",
      "target_type": "identity"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Education Sector Organizations",
      "target_type": "identity"
    },
    {
      "source_name": "KongTuke",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "IT Sector Organizations",
      "target_type": "identity"
    },
    {
      "source_name": "Mistic",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "DLL Side-Loading",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mistic",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Match Legitimate Resource Name or Location",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mistic",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "File Deletion",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mistic",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Web Protocols",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mistic",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "Mistic C2 Infrastructure",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Qilin",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Qilin",
      "target_type": "malware"
    },
    {
      "source_name": "Akira",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Akira",
      "target_type": "malware"
    },
    {
      "source_name": "Black Basta",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Black Basta",
      "target_type": "malware"
    },
    {
      "source_name": "Mistic Backdoor Persistence via Run Key",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Mistic",
      "target_type": "malware"
    },
    {
      "source_name": "Mistic Backdoor Persistence via Run Key",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Mistic Backdoor Persistence via Run Key",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "target_type": "windows-registry-key"
    },
    {
      "source_name": "Mistic Backdoor Associated File Names",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Mistic",
      "target_type": "malware"
    },
    {
      "source_name": "Mistic Backdoor Associated File Names",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "EndpointDlp.dll",
      "target_type": "file"
    },
    {
      "source_name": "Mistic Backdoor Associated File Names",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "MpExtMs.exe",
      "target_type": "file"
    }
  ]
}