# stix/sidecopy-apt-targets-afghanistans-finance-ministry-with-xenorat-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-06-02T15:21:39.150Z",
    "article_slug": "sidecopy-apt-targets-afghanistans-finance-ministry-with-xenorat",
    "article_date": "2026-06-02",
    "token_usage": {
      "input_tokens": 144172,
      "output_tokens": 3866,
      "cached_tokens": 137782,
      "total_tokens": 148038,
      "cache_hit": true,
      "estimated_cost_usd": 0.06387
    }
  },
  "report_name": "SideCopy APT Targets Afghanistan's Finance Ministry in 'XENOFISCAL' Espionage Campaign",
  "threat_actors": [
    {
      "name": "SideCopy",
      "mitre_id": "G1008",
      "aliases": [],
      "description": "A Pakistan-aligned APT group, considered a subgroup of Transparent Tribe (APT36), conducting the 'Operation XENOFISCAL' cyber-espionage campaign against Afghanistan's Ministry of Finance. The group uses spear-phishing with Pashto-language lures to deploy XenoRAT.",
      "goals": "Cyber-espionage, data theft, and surveillance against government and military targets in South Asia.",
      "primary_motivation": "organizational-gain",
      "sophistication": "intermediate"
    },
    {
      "name": "Transparent Tribe",
      "mitre_id": "G0134",
      "aliases": [
        "APT36"
      ],
      "description": "A Pakistan-aligned threat group, also known as APT36, that acts as an umbrella group for other actors, including SideCopy. This campaign is attributed to its SideCopy subgroup."
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Attachment",
      "mitre_id": "T1566.001",
      "description": "The attack begins with a spear-phishing email containing a ZIP archive with a malicious LNK file, sent to officials at Afghanistan's Ministry of Finance."
    },
    {
      "name": "Malicious File",
      "mitre_id": "T1204.002",
      "description": "Victims are tricked into executing a malicious LNK file disguised with a Pashto filename, which initiates the infection chain."
    },
    {
      "name": "System Binary Proxy Execution: Mshta",
      "mitre_id": "T1218.005",
      "description": "The malicious LNK file executes mshta.exe to download and run a remote HTML Application (HTA) file from a compromised domain."
    },
    {
      "name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "mitre_id": "T1547.001",
      "description": "XenoRAT establishes persistence by creating a registry run key, often mimicking a legitimate application like Microsoft Edge, to ensure it executes on system startup."
    },
    {
      "name": "Application Layer Protocol: Web Protocols",
      "mitre_id": "T1071.001",
      "description": "The infection chain uses mshta.exe to download a remote HTA file over HTTP. The final XenoRAT payload communicates with its C2 server over TCP, which may use web protocols."
    },
    {
      "name": "Obtain Capabilities: Tool",
      "mitre_id": "T1588.002",
      "description": "SideCopy obtained the open-source XenoRAT for use as the final payload in their espionage campaign."
    }
  ],
  "malware": [
    {
      "name": "XenoRAT",
      "mitre_id": "",
      "type": "remote-access-trojan",
      "is_family": true,
      "aliases": [],
      "description": "An open-source remote access trojan (version 1.8.7 mentioned) deployed as the final payload. It provides attackers with full remote access for keylogging, screen capture, file system access, and data exfiltration.",
      "capabilities": [
        "communicates-with-c2",
        "captures-input-peripherals",
        "captures-output-peripherals",
        "controls-local-machine",
        "exfiltrates-data",
        "persists-after-system-reboot"
      ]
    }
  ],
  "tools": [
    {
      "name": "mshta.exe",
      "mitre_id": "S0196",
      "aliases": [],
      "description": "The legitimate Windows utility mshta.exe was abused by the attackers to execute a remote HTML Application (HTA) file from a compromised server as part of the infection chain.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [
    {
      "name": "Operation XENOFISCAL",
      "description": "A targeted cyber-espionage campaign by the SideCopy APT group against Afghanistan's Ministry of Finance, using spear-phishing to deploy the XenoRAT trojan.",
      "objective": "Espionage and data exfiltration from the Afghan Ministry of Finance."
    }
  ],
  "identities": [
    {
      "name": "Afghanistan Ministry of Finance",
      "identity_class": "organization",
      "sectors": [
        "government"
      ],
      "description": "The primary target of the 'Operation XENOFISCAL' campaign. The attack specifically focused on its 34 provincial revenue and finance directorates (Mustoufiats) and Pashto-speaking officials."
    }
  ],
  "locations": [
    {
      "name": "Afghanistan",
      "region": "southern-asia"
    },
    {
      "name": "Pakistan",
      "region": "southern-asia"
    },
    {
      "name": "South Asia",
      "region": "southern-asia"
    }
  ],
  "infrastructure": [
    {
      "name": "Malicious HTA Hosting Server",
      "infrastructure_types": [
        "hosting-malware",
        "command-and-control"
      ],
      "description": "A compromised Afghan educational domain (abimj.edu.af) used to host the malicious HTA file for the initial stage of the infection."
    },
    {
      "name": "XenoRAT C2 Server",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "Command and control server (185.235.137.106) used by the XenoRAT payload to receive commands and exfiltrate data."
    }
  ],
  "observables": [
    {
      "type": "domain-name",
      "value": "abimj.edu.af"
    },
    {
      "type": "ipv4-addr",
      "value": "185.235.137.106"
    },
    {
      "type": "windows-registry-key",
      "value": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
    }
  ],
  "indicators": [
    {
      "name": "SideCopy XENOFISCAL C2 and Staging Infrastructure",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[domain-name:value = 'abimj.edu.af'] OR [ipv4-addr:value = '185.235.137.106']",
      "x_netsecio_sigma_rule": "title: SideCopy XENOFISCAL C2 and Staging Network Activity\nid: 4f5e6d7c-8b9a-4123-b456-7890abcdef12\nstatus: experimental\ndescription: Detects network connections to domains and IPs associated with the SideCopy 'Operation XENOFISCAL' campaign.\nreferences:\n  - https://thehackernews.com/2026/06/pakistan-linked-sidecopy-targets.html\nauthor: cyber.netsecops.io\ndate: 2026-06-02\ntags:\n  - attack.t1071.001\n  - attack.g1008\nlevel: high\nlogsource:\n  category: network_traffic\n  product: windows\ndetection:\n  selection_domain:\n    destination.domain: 'abimj.edu.af'\n  selection_ip:\n    destination.ip: '185.235.137.106'\n  condition: 1 of them",
      "pattern": "",
      "description": "Detects network traffic to the compromised HTA staging domain or the XenoRAT C2 server IP address identified in the Operation XENOFISCAL campaign."
    },
    {
      "name": "Suspicious Remote HTA Execution via Mshta",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Remote HTA Execution via Mshta.exe\nid: 9a8b7c6d-5e4f-4321-a0b9-c8d7e6f5a4b3\nstatus: experimental\ndescription: Detects the use of mshta.exe to execute a remote HTML Application (HTA) file, a technique used by SideCopy in the XENOFISCAL campaign.\nreferences:\n  - https://thehackernews.com/2026/06/pakistan-linked-sidecopy-targets.html\nauthor: cyber.netsecops.io\ndate: 2026-06-02\ntags:\n  - attack.t1218.005\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\mshta.exe'\n    CommandLine|contains:\n      - 'http://'\n      - 'https://'\n  condition: selection",
      "description": "Detects the execution of mshta.exe with a command line argument containing a URL, which is indicative of remote HTA file execution used in this campaign."
    },
    {
      "name": "Persistence via Registry Run Key",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Persistence via HKCU Run Key Creation\nid: 1a2b3c4d-5e6f-4789-890a-bcdef0123456\nstatus: experimental\ndescription: 'Detects the creation of a new value in the HKCU\\...\\Run registry key, a common persistence technique used by malware like XenoRAT.'\nreferences:\n  - https://thehackernews.com/2026/06/pakistan-linked-sidecopy-targets.html\nauthor: cyber.netsecops.io\ndate: 2026-06-02\ntags:\n  - attack.t1547.001\nlevel: medium\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    EventType: 'SetValue'\n    TargetObject|contains: '\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n  condition: selection",
      "description": "Detects the creation of a new persistence entry in the HKEY_CURRENT_USER Run key, a method used by XenoRAT in this campaign."
    }
  ],
  "relationships": [
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "attributed-to",
      "target_name": "Transparent Tribe",
      "target_type": "threat-actor"
    },
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "XenoRAT",
      "target_type": "malware"
    },
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "mshta.exe",
      "target_type": "tool"
    },
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Afghanistan Ministry of Finance",
      "target_type": "identity"
    },
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "originates-from",
      "target_name": "Pakistan",
      "target_type": "location"
    },
    {
      "source_name": "SideCopy",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "South Asia",
      "target_type": "location"
    },
    {
      "source_name": "Operation XENOFISCAL",
      "source_type": "campaign",
      "relationship_type": "attributed-to",
      "target_name": "SideCopy",
      "target_type": "threat-actor"
    },
    {
      "source_name": "Operation XENOFISCAL",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "Afghanistan Ministry of Finance",
      "target_type": "identity"
    },
    {
      "source_name": "Operation XENOFISCAL",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Spearphishing Attachment",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "XenoRAT",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "XenoRAT",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "XenoRAT C2 Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Afghanistan Ministry of Finance",
      "source_type": "identity",
      "relationship_type": "located-at",
      "target_name": "Afghanistan",
      "target_type": "location"
    },
    {
      "source_name": "SideCopy XENOFISCAL C2 and Staging Infrastructure",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "SideCopy",
      "target_type": "threat-actor"
    },
    {
      "source_name": "SideCopy XENOFISCAL C2 and Staging Infrastructure",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "abimj.edu.af",
      "target_type": "domain-name"
    },
    {
      "source_name": "SideCopy XENOFISCAL C2 and Staging Infrastructure",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "185.235.137.106",
      "target_type": "ipv4-addr"
    },
    {
      "source_name": "Suspicious Remote HTA Execution via Mshta",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "System Binary Proxy Execution: Mshta",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Persistence via Registry Run Key",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    }
  ]
}