# stix/sleepergem-supply-chain-attack-hijacks-dormant-rubygems-accounts-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-07-20T15:42:49.072Z",
    "article_slug": "sleepergem-supply-chain-attack-hijacks-dormant-rubygems-accounts",
    "article_date": "2026-07-20",
    "token_usage": {
      "input_tokens": 146653,
      "output_tokens": 4172,
      "cached_tokens": 140179,
      "total_tokens": 150825,
      "cache_hit": true,
      "estimated_cost_usd": 0.067335
    }
  },
  "report_name": "'SleeperGem' Supply Chain Attack Hits RubyGems via Dormant Accounts",
  "threat_actors": [],
  "attack_patterns": [
    {
      "name": "Compromise Software Dependencies and Development Tools",
      "mitre_id": "T1195.001",
      "description": "Attackers published malicious versions of legitimate RubyGems packages by adding a malicious dependency, `git_credential_manager`."
    },
    {
      "name": "Valid Accounts",
      "mitre_id": "T1078",
      "description": "Attackers hijacked long-dormant maintainer accounts (e.g., 'LR-DEV', 'pinkroom') on RubyGems to publish malicious package versions."
    },
    {
      "name": "System Checks",
      "mitre_id": "T1497.001",
      "description": "The malicious code included checks to determine if it was running in a CI/CD or automated build environment, aborting execution to evade automated security scanning and specifically target developer machines."
    },
    {
      "name": "Systemd Service",
      "mitre_id": "T1543.002",
      "description": "The second-stage backdoor payload established persistence on developer machines by creating malicious systemd services."
    },
    {
      "name": "Cron",
      "mitre_id": "T1053.003",
      "description": "The second-stage backdoor payload established persistence on developer machines by creating malicious cron jobs."
    },
    {
      "name": "Setuid and Setgid",
      "mitre_id": "T1548.001",
      "description": "The backdoor attempted to escalate privileges to root by installing a setuid root shell at /usr/local/sbin/ping6."
    }
  ],
  "malware": [
    {
      "name": "git_credential_manager",
      "mitre_id": "",
      "type": "dropper",
      "is_family": false,
      "aliases": [],
      "description": "A malicious Ruby gem (versions 2.8.0-2.8.3) impersonating Microsoft's Git Credential Manager. It acted as a first-stage dropper to deliver a persistent backdoor and was added as a dependency to other compromised gems.",
      "capabilities": [
        "installs-other-components",
        "anti-sandbox"
      ]
    },
    {
      "name": "SleeperGem Backdoor",
      "mitre_id": "",
      "type": "backdoor",
      "is_family": false,
      "aliases": [],
      "description": "A persistent, multi-stage backdoor installed by the git_credential_manager dropper. It establishes persistence via systemd and cron, and attempts to escalate privileges by creating a setuid root shell.",
      "capabilities": [
        "persists-after-system-reboot",
        "escalates-privileges"
      ]
    }
  ],
  "tools": [
    {
      "name": "Microsoft Git Credential Manager",
      "mitre_id": "",
      "aliases": [],
      "description": "A legitimate Microsoft tool that was impersonated by the malicious `git_credential_manager` Ruby gem as part of a dependency confusion attack."
    },
    {
      "name": "Dendreo",
      "mitre_id": "",
      "aliases": [],
      "description": "A legitimate Ruby gem that was compromised in versions 1.1.3 and 1.1.4 to include the malicious `git_credential_manager` gem as a dependency, thereby distributing the malware."
    },
    {
      "name": "fastlane-plugin-run_tests_firebase_testlab",
      "mitre_id": "",
      "aliases": [],
      "description": "A popular, legitimate Ruby gem that was compromised in version 0.3.2 to include the malicious `git_credential_manager` gem as a dependency, significantly amplifying the attack's reach to over 574,000 downloads."
    }
  ],
  "vulnerabilities": [],
  "campaigns": [
    {
      "name": "SleeperGem",
      "description": "A sophisticated software supply chain attack targeting the RubyGems ecosystem by hijacking dormant maintainer accounts to publish malicious versions of legitimate packages and distribute a persistent backdoor onto developer workstations.",
      "first_seen": "2026-07-18T00:00:00Z"
    }
  ],
  "identities": [
    {
      "name": "RubyGems",
      "identity_class": "organization",
      "sectors": [
        "technology"
      ],
      "description": "The targeted software package registry for the Ruby programming language, which was compromised through the takeover of dormant maintainer accounts."
    },
    {
      "name": "Developers",
      "identity_class": "class",
      "sectors": [
        "technology"
      ],
      "description": "The intended victims of the SleeperGem attack, whose workstations were targeted for backdoor installation, credential theft, and potential lateral movement."
    },
    {
      "name": "LR-DEV",
      "identity_class": "individual",
      "sectors": [
        "technology"
      ],
      "description": "A dormant RubyGems maintainer account that was hijacked by attackers to publish malicious packages as part of the SleeperGem campaign."
    },
    {
      "name": "pinkroom",
      "identity_class": "individual",
      "sectors": [
        "technology"
      ],
      "description": "A dormant RubyGems maintainer account that was hijacked by attackers to publish malicious packages as part of the SleeperGem campaign."
    }
  ],
  "locations": [],
  "infrastructure": [
    {
      "name": "SleeperGem Forgejo C2 Host",
      "infrastructure_types": [
        "command-and-control",
        "hosting-malware"
      ],
      "description": "An attacker-controlled Forgejo (Git service) host used by early versions of the git_credential_manager dropper to fetch a second-stage payload."
    }
  ],
  "observables": [
    {
      "type": "file",
      "name": "git_credential_manager"
    },
    {
      "type": "file",
      "name": "Dendreo"
    },
    {
      "type": "file",
      "name": "fastlane-plugin-run_tests_firebase_testlab"
    },
    {
      "type": "file",
      "name": "/usr/local/sbin/ping6"
    }
  ],
  "indicators": [
    {
      "name": "SleeperGem Malicious RubyGem Installation",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[process:command_line MATCHES 'gem install git_credential_manager']",
      "x_netsecio_sigma_rule": "title: SleeperGem Malicious RubyGem Installation\nid: f9a8b7c6-d5e4-4f3a-2b1c-0a9b8c7d6e5f\nstatus: experimental\ndescription: Detects the installation of the malicious 'git_credential_manager' Ruby gem, a component of the SleeperGem supply chain attack.\nreferences:\n  - https://thehackernews.com/2026/07/sleepergem-uses-three-malicious.html\nauthor: cyber.netsecops.io\ndate: 2026-07-20\ntags:\n  - attack.t1195.001\n  - attack.initial_access\nlevel: high\nlogsource:\n  category: process_creation\n  product: linux\ndetection:\n  selection:\n    Image|endswith: '/gem'\n    CommandLine|contains|all:\n      - 'install'\n      - 'git_credential_manager'\n  condition: selection",
      "pattern": ""
    },
    {
      "name": "SleeperGem Setuid Root Shell Creation",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: SleeperGem Setuid Root Shell Creation\nid: e8d7c6b5-a4f3-4e2d-1c0b-9a8b7c6d5e4f\nstatus: experimental\ndescription: Detects the use of chmod to set the SUID bit on /usr/local/sbin/ping6, a technique used by the SleeperGem backdoor for privilege escalation.\nreferences:\n  - https://thehackernews.com/2026/07/sleepergem-uses-three-malicious.html\nauthor: cyber.netsecops.io\ndate: 2026-07-20\ntags:\n  - attack.t1548.001\n  - attack.privilege_escalation\nlevel: high\nlogsource:\n  category: process_creation\n  product: linux\ndetection:\n  selection:\n    Image|endswith: '/chmod'\n    CommandLine|contains|all:\n      - 'u+s'\n      - '/usr/local/sbin/ping6'\n  condition: selection"
    },
    {
      "name": "SleeperGem Systemd User Service Persistence",
      "indicator_types": [
        "malicious-activity",
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: SleeperGem Systemd User Service Persistence\nid: d7c6b5a4-93e2-4d1c-0b9a-8b7c6d5e4f3a\nstatus: experimental\ndescription: Detects the creation of a new service file in a user's systemd directory, a persistence technique used by the SleeperGem backdoor.\nreferences:\n  - https://thehackernews.com/2026/07/sleepergem-uses-three-malicious.html\nauthor: cyber.netsecops.io\ndate: 2026-07-20\ntags:\n  - attack.t1543.002\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: file_event\n  product: linux\ndetection:\n  selection:\n    TargetFilename|contains: '/.config/systemd/user/'\n    TargetFilename|endswith: '.service'\n  condition: selection"
    },
    {
      "name": "SleeperGem Cron Job Persistence",
      "indicator_types": [
        "malicious-activity",
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: SleeperGem Cron Job Persistence\nid: c6b5a493-82d1-4c0b-9a8b-7c6d5e4f3a2b\nstatus: experimental\ndescription: Detects suspicious modification of user crontab files, a persistence technique used by the SleeperGem backdoor.\nreferences:\n  - https://thehackernews.com/2026/07/sleepergem-uses-three-malicious.html\nauthor: cyber.netsecops.io\ndate: 2026-07-20\ntags:\n  - attack.t1053.003\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: process_creation\n  product: linux\ndetection:\n  selection:\n    Image|endswith: '/crontab'\n  filter_legit:\n    CommandLine|contains:\n      - ' -l'\n      - ' -r'\n  condition: selection and not filter_legit"
    }
  ],
  "relationships": [
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "git_credential_manager",
      "target_type": "malware"
    },
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "SleeperGem Backdoor",
      "target_type": "malware"
    },
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "RubyGems",
      "target_type": "identity"
    },
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "targets",
      "target_name": "Developers",
      "target_type": "identity"
    },
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Compromise Software Dependencies and Development Tools",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Valid Accounts",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "git_credential_manager",
      "source_type": "malware",
      "relationship_type": "impersonates",
      "target_name": "Microsoft Git Credential Manager",
      "target_type": "tool"
    },
    {
      "source_name": "git_credential_manager",
      "source_type": "malware",
      "relationship_type": "downloads",
      "target_name": "SleeperGem Backdoor",
      "target_type": "malware"
    },
    {
      "source_name": "git_credential_manager",
      "source_type": "malware",
      "relationship_type": "downloads-from",
      "target_name": "SleeperGem Forgejo C2 Host",
      "target_type": "infrastructure"
    },
    {
      "source_name": "git_credential_manager",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "System Checks",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Dendreo",
      "source_type": "tool",
      "relationship_type": "delivers",
      "target_name": "git_credential_manager",
      "target_type": "malware"
    },
    {
      "source_name": "fastlane-plugin-run_tests_firebase_testlab",
      "source_type": "tool",
      "relationship_type": "delivers",
      "target_name": "git_credential_manager",
      "target_type": "malware"
    },
    {
      "source_name": "SleeperGem Backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Systemd Service",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem Backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Cron",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem Backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Setuid and Setgid",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem Malicious RubyGem Installation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "git_credential_manager",
      "target_type": "malware"
    },
    {
      "source_name": "SleeperGem Setuid Root Shell Creation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Setuid and Setgid",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem Systemd User Service Persistence",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Systemd Service",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "SleeperGem Cron Job Persistence",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Cron",
      "target_type": "attack-pattern"
    }
  ]
}