# stix/social-engineering-campaign-abuses-whatsapp-for-windows-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T16:09:49.567Z",
    "article_slug": "social-engineering-campaign-abuses-whatsapp-for-windows",
    "article_date": "2026-04-02",
    "token_usage": {
      "input_tokens": 141544,
      "output_tokens": 3982,
      "cached_tokens": 135916,
      "total_tokens": 145526,
      "cache_hit": true,
      "estimated_cost_usd": 0.063844
    }
  },
  "report_name": "Microsoft Warns of Social Engineering Campaign Abusing WhatsApp for Windows",
  "threat_actors": [
    {
      "name": "Unnamed Social Engineering Actor",
      "mitre_id": "",
      "aliases": [],
      "description": "An unnamed threat actor conducting a social engineering campaign targeting WhatsApp for Windows users. The actor uses VBScript malware to install remote access software, relying on living-off-the-land techniques.",
      "sophistication": "intermediate"
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Attachment",
      "mitre_id": "T1566.001",
      "description": "The initial delivery mechanism is a malicious VBScript attachment sent via the WhatsApp messaging service, disguised as a legitimate document or image."
    },
    {
      "name": "User Execution: Malicious File",
      "mitre_id": "T1204.002",
      "description": "The attack's success is entirely dependent on the user being tricked into executing the initial malicious .vbs file."
    },
    {
      "name": "Command and Scripting Interpreter: Visual Basic",
      "mitre_id": "T1059.005",
      "description": "The initial payload is a Visual Basic Script (.vbs) file executed by the Windows Script Host (wscript.exe)."
    },
    {
      "name": "System Binary Proxy Execution",
      "mitre_id": "T1218",
      "description": "The malware uses renamed, legitimate Windows binaries such as bitsadmin or certutil to download a second-stage payload from an attacker-controlled server, evading simple detection rules."
    },
    {
      "name": "Abuse Elevation Control Mechanism: Bypass User Account Control",
      "mitre_id": "T1548.002",
      "description": "The malware modifies User Account Control (UAC) settings, such as the EnableLUA registry value, to operate with fewer restrictions and suppress future warnings."
    },
    {
      "name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "mitre_id": "T1547.001",
      "description": "The malware creates registry entries in autostart locations, such as the Run key, to achieve persistence on the compromised machine."
    }
  ],
  "malware": [
    {
      "name": "VBScript Downloader",
      "mitre_id": "",
      "type": "downloader",
      "is_family": false,
      "aliases": [],
      "description": "A malicious Visual Basic Script (.vbs) file delivered via WhatsApp. When executed, it uses living-off-the-land techniques to download and run a second-stage remote access payload.",
      "capabilities": [
        "installs-other-components",
        "communicates-with-c2"
      ]
    },
    {
      "name": "Remote Access Software",
      "mitre_id": "",
      "type": "remote-access-trojan",
      "is_family": false,
      "aliases": [],
      "description": "The final payload installed by the VBScript downloader. It is delivered as an unsigned Microsoft Installer (MSI) package and gives the attacker full, persistent access to the victim's system.",
      "capabilities": [
        "controls-local-machine",
        "escalates-privileges",
        "persists-after-system-reboot"
      ]
    }
  ],
  "tools": [
    {
      "name": "wscript.exe",
      "mitre_id": "",
      "aliases": [],
      "description": "The Windows Script Host (wscript.exe) is used to execute the initial malicious .vbs file.",
      "tool_types": []
    },
    {
      "name": "bitsadmin",
      "mitre_id": "S0190",
      "aliases": [],
      "description": "The legitimate Windows BITSAdmin tool is copied, renamed, and used by the VBScript malware to download a second-stage payload from an attacker-controlled server.",
      "tool_types": []
    },
    {
      "name": "certutil",
      "mitre_id": "S0160",
      "aliases": [],
      "description": "The legitimate Windows certutil tool is copied, renamed, and potentially used by the VBScript malware to download a second-stage payload.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [],
  "identities": [],
  "locations": [],
  "infrastructure": [
    {
      "name": "Attacker-controlled Server",
      "infrastructure_types": [
        "command-and-control",
        "hosting-malware"
      ],
      "description": "A remote server controlled by the attacker used to host and deliver the second-stage malware payload."
    }
  ],
  "observables": [
    {
      "type": "windows-registry-key",
      "value": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
    },
    {
      "type": "windows-registry-key",
      "value": "HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableLUA"
    }
  ],
  "indicators": [
    {
      "name": "Suspicious VBScript Execution by WScript",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious VBScript Execution by WScript\nid: 2b8f1c9d-5a3e-4b1f-8c9a-1d2e3f4a5b6c\nstatus: experimental\ndescription: Detects the execution of wscript.exe with a .vbs file argument, which is a common technique for initial execution in social engineering campaigns.\nreferences:\n  - https://www.malwarebytes.com/blog/news/2026/04/whatsapp-on-windows-users-targeted-in-new-campaign-warns-microsoft\nauthor: cyber.netsecops.io\ndate: 2026-04-02\ntags:\n  - attack.t1059.005\n  - attack.execution\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\wscript.exe'\n    CommandLine|contains: '.vbs'\n  condition: selection",
      "description": "Detects the execution of wscript.exe with a .vbs file, as seen in the WhatsApp social engineering campaign."
    },
    {
      "name": "File Download via BITSAdmin",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: File Download via BITSAdmin\nid: 8c9a1d2e-3f4a-4b6c-9d2b-8f1c9d5a3e4b\nstatus: experimental\ndescription: Detects the use of bitsadmin.exe with the /transfer switch to download a remote file, a common living-off-the-land technique.\nreferences:\n  - https://www.malwarebytes.com/blog/news/2026/04/whatsapp-on-windows-users-targeted-in-new-campaign-warns-microsoft\nauthor: cyber.netsecops.io\ndate: 2026-04-02\ntags:\n  - attack.t1218\n  - attack.defense_evasion\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\bitsadmin.exe'\n    CommandLine|contains:\n      - ' /transfer '\n      - ' /addfile '\n  condition: selection",
      "description": "Detects the use of bitsadmin.exe to download payloads, a technique used in this campaign."
    },
    {
      "name": "UAC Bypass via EnableLUA Registry Modification",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[windows-registry-key:key = 'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Policies\\\\System\\\\EnableLUA']",
      "x_netsecio_sigma_rule": "title: UAC Bypass via EnableLUA Registry Modification\nid: 5a3e4b1f-8c9a-4d2e-b8f1-c9d5a3e4b6c9\nstatus: experimental\ndescription: Monitors for modifications to the EnableLUA registry key, which can be used to disable User Account Control (UAC).\nreferences:\n  - https://www.malwarebytes.com/blog/news/2026/04/whatsapp-on-windows-users-targeted-in-new-campaign-warns-microsoft\nauthor: cyber.netsecops.io\ndate: 2026-04-02\ntags:\n  - attack.t1548.002\n  - attack.privilege_escalation\nlevel: high\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    TargetObject|endswith: '\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableLUA'\n    Details: '0'\n  condition: selection",
      "pattern": "",
      "description": "Detects attempts to disable UAC by modifying the EnableLUA registry value, a technique used by the malware in this campaign."
    },
    {
      "name": "Persistence via Registry Run Key",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "[windows-registry-key:key MATCHES '^HKEY_CURRENT_USER\\\\\\\\Software\\\\\\\\Microsoft\\\\\\\\Windows\\\\\\\\CurrentVersion\\\\\\\\Run\\\\\\\\.*']",
      "x_netsecio_sigma_rule": "title: Persistence via Registry Run Key Creation\nid: 1d2e3f4a-5b6c-4b8f-8c9a-1d2e3f4a5b6d\nstatus: experimental\ndescription: Detects the creation of a new value in the HKCU Run registry key, a common persistence technique.\nreferences:\n  - https://www.malwarebytes.com/blog/news/2026/04/whatsapp-on-windows-users-targeted-in-new-campaign-warns-microsoft\nauthor: cyber.netsecops.io\ndate: 2026-04-02\ntags:\n  - attack.t1547.001\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    EventType: 'SetValue'\n    TargetObject|startswith: 'HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n  condition: selection",
      "pattern": "",
      "description": "Detects the creation of a new entry in the user's registry Run key for persistence, as used in this campaign."
    }
  ],
  "relationships": [
    {
      "source_name": "Unnamed Social Engineering Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Spearphishing Attachment",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Unnamed Social Engineering Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "User Execution: Malicious File",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Unnamed Social Engineering Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "VBScript Downloader",
      "target_type": "malware"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Command and Scripting Interpreter: Visual Basic",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "System Binary Proxy Execution",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "bitsadmin",
      "target_type": "tool"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "certutil",
      "target_type": "tool"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "downloads",
      "target_name": "Remote Access Software",
      "target_type": "malware"
    },
    {
      "source_name": "VBScript Downloader",
      "source_type": "malware",
      "relationship_type": "communicates-with",
      "target_name": "Attacker-controlled Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Remote Access Software",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Abuse Elevation Control Mechanism: Bypass User Account Control",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Remote Access Software",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Suspicious VBScript Execution by WScript",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Command and Scripting Interpreter: Visual Basic",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "File Download via BITSAdmin",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "System Binary Proxy Execution",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC Bypass via EnableLUA Registry Modification",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Abuse Elevation Control Mechanism: Bypass User Account Control",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Persistence via Registry Run Key",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC Bypass via EnableLUA Registry Modification",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableLUA",
      "target_type": "windows-registry-key"
    },
    {
      "source_name": "Persistence via Registry Run Key",
      "source_type": "indicator",
      "relationship_type": "based-on",
      "target_name": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "target_type": "windows-registry-key"
    }
  ]
}