# stix/stealthy-powmix-botnet-targets-czech-republic-workforce-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T16:59:58.053Z",
    "article_slug": "stealthy-powmix-botnet-targets-czech-republic-workforce",
    "article_date": "2026-04-17",
    "token_usage": {
      "input_tokens": 141797,
      "output_tokens": 3242,
      "cached_tokens": 136033,
      "total_tokens": 145039,
      "cache_hit": true,
      "estimated_cost_usd": 0.056629
    }
  },
  "report_name": "Stealthy 'PowMix' Botnet Targets Czech Workforce with Evasive C2 Communications",
  "threat_actors": [
    {
      "name": "PowMix Actor",
      "mitre_id": "",
      "aliases": [],
      "description": "The unspecified threat actor behind the PowMix botnet campaign, targeting professionals in the Czech Republic since at least December 2025. The actor demonstrates sophistication in social engineering and technical implementation, with a focus on stealthy C2 communications to evade detection.",
      "sophistication": "intermediate"
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Attachment",
      "mitre_id": "T1566.001",
      "description": "The attack begins with a phishing email containing a malicious ZIP archive with a .LNK file, using compliance-themed lures to target professionals."
    },
    {
      "name": "Malicious File",
      "mitre_id": "T1204.002",
      "description": "Victims are tricked into opening a malicious Windows Shortcut (.LNK) file contained within a ZIP attachment, which initiates the execution chain."
    },
    {
      "name": "Command and Scripting Interpreter: PowerShell",
      "mitre_id": "T1059.001",
      "description": "A PowerShell loader is executed by the malicious LNK file to decrypt and run the final PowMix payload directly in memory, evading file-based detection."
    },
    {
      "name": "Scheduled Task/Job: Scheduled Task",
      "mitre_id": "T1053.005",
      "description": "The PowMix malware establishes persistence on compromised systems by creating a scheduled task to ensure it runs again after a reboot."
    },
    {
      "name": "Application Layer Protocol: Web Protocols",
      "mitre_id": "T1071.001",
      "description": "PowMix uses HTTP/S for C2 communications, embedding encrypted data in URL paths to mimic legitimate REST API calls and evade network signature detection."
    }
  ],
  "malware": [
    {
      "name": "PowMix",
      "mitre_id": "",
      "type": "botnet",
      "is_family": true,
      "aliases": [],
      "description": "A botnet delivered via phishing that uses PowerShell for in-memory execution. It is designed for stealth, employing randomized C2 beaconing intervals and embedding encrypted data into URL paths to mimic legitimate API traffic.",
      "capabilities": [
        "communicates-with-c2",
        "persists-after-system-reboot"
      ]
    }
  ],
  "tools": [
    {
      "name": "PowerShell",
      "mitre_id": "S0194",
      "aliases": [],
      "description": "Used as a loader to decrypt and run the final PowMix malware payload directly in memory.",
      "tool_types": []
    },
    {
      "name": "schtasks",
      "mitre_id": "S0111",
      "aliases": [
        "schtasks.exe"
      ],
      "description": "Used by the PowMix malware to create a scheduled task for persistence on compromised systems.",
      "tool_types": []
    }
  ],
  "vulnerabilities": [],
  "campaigns": [
    {
      "name": "PowMix Campaign",
      "description": "An ongoing campaign active since at least December 2025, targeting the workforce in the Czech Republic with the PowMix botnet via phishing emails.",
      "aliases": [],
      "first_seen": "2025-12-01T00:00:00Z"
    }
  ],
  "identities": [
    {
      "name": "Czech Professionals in HR, Legal, and Finance",
      "identity_class": "class",
      "sectors": [
        "professional-services",
        "legal",
        "financial-services"
      ],
      "description": "The primary target group for the PowMix botnet campaign, lured with compliance-themed phishing emails."
    },
    {
      "name": "EDEKA",
      "identity_class": "organization",
      "sectors": [
        "retail"
      ],
      "description": "A legitimate brand impersonated in phishing lures to add credibility to the campaign targeting professionals in the Czech Republic."
    }
  ],
  "locations": [
    {
      "name": "Czech Republic",
      "country": "Czech Republic",
      "region": "eastern-europe"
    }
  ],
  "infrastructure": [
    {
      "name": "PowMix C2 Infrastructure",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "Command and control servers used by the PowMix botnet. Communications are designed for stealth, using randomized beaconing intervals and embedding encrypted data in URL paths."
    }
  ],
  "observables": [],
  "indicators": [
    {
      "name": "PowMix PowerShell Loader Execution",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: PowMix PowerShell Loader Execution\nid: 9c5b1a3e-f8d2-4e1c-b9a5-6f7d8e9c0b1a\nstatus: experimental\ndescription: Detects the execution of a PowerShell script with a bypassed execution policy, a technique used by PowMix's LNK-based loader.\nreferences:\n  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html\nauthor: cyber.netsecops.io\ndate: 2026-04-17\ntags:\n  - attack.t1059.001\n  - attack.execution\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\powershell.exe'\n    CommandLine|contains|all:\n      - '-ExecutionPolicy Bypass'\n      - '-File'\n  condition: selection",
      "description": "Detects the execution of a PowerShell script with a bypassed execution policy, a technique used by PowMix's LNK-based loader."
    },
    {
      "name": "PowMix Persistence via Scheduled Task",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: PowMix Persistence via Scheduled Task\nid: 3d2e1c0b-a9b8-4f7d-8e6c-5b4a3d2e1c0b\nstatus: experimental\ndescription: Detects the creation of a new scheduled task using schtasks.exe, which is the persistence mechanism used by the PowMix botnet.\nreferences:\n  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html\nauthor: cyber.netsecops.io\ndate: 2026-04-17\ntags:\n  - attack.t1053.005\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\schtasks.exe'\n    CommandLine|contains: '/Create'\n  condition: selection",
      "description": "Detects the creation of a new scheduled task using schtasks.exe, which is the persistence mechanism used by the PowMix botnet."
    },
    {
      "name": "PowMix C2 URL Pattern",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: PowMix C2 URL Pattern\nid: 7a6b5c4d-3e2f-4d1c-a9b8-c7d6e5f4a3b2\nstatus: experimental\ndescription: Detects outbound web requests with URL paths matching the pattern used by PowMix for C2 communications, which embeds data to mimic REST API calls.\nreferences:\n  - https://thehackernews.com/2026/04/newly-discovered-powmix-botnet-hits.html\nauthor: cyber.netsecops.io\ndate: 2026-04-17\ntags:\n  - attack.t1071.001\n  - attack.command_and_control\nlevel: medium\nlogsource:\n  category: proxy\ndetection:\n  selection:\n    c-uri|re: '/api/v\\d+/data=.*'\n  condition: selection",
      "description": "Detects outbound web requests with URL paths matching the pattern used by PowMix for C2 communications, which embeds data to mimic REST API calls."
    }
  ],
  "relationships": [
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "PowMix",
      "target_type": "malware"
    },
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "PowerShell",
      "target_type": "tool"
    },
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "schtasks",
      "target_type": "tool"
    },
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Czech Professionals in HR, Legal, and Finance",
      "target_type": "identity"
    },
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Czech Republic",
      "target_type": "location"
    },
    {
      "source_name": "PowMix Actor",
      "source_type": "threat-actor",
      "relationship_type": "impersonates",
      "target_name": "EDEKA",
      "target_type": "identity"
    },
    {
      "source_name": "PowMix Campaign",
      "source_type": "campaign",
      "relationship_type": "attributed-to",
      "target_name": "PowMix Actor",
      "target_type": "threat-actor"
    },
    {
      "source_name": "PowMix Campaign",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Spearphishing Attachment",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "PowMix Campaign",
      "source_type": "campaign",
      "relationship_type": "uses",
      "target_name": "Malicious File",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "PowMix",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Command and Scripting Interpreter: PowerShell",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "PowMix",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Scheduled Task/Job: Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "PowMix",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Application Layer Protocol: Web Protocols",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "PowMix",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "PowMix C2 Infrastructure",
      "target_type": "infrastructure"
    },
    {
      "source_name": "PowMix PowerShell Loader Execution",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "PowMix",
      "target_type": "malware"
    },
    {
      "source_name": "PowMix Persistence via Scheduled Task",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "PowMix",
      "target_type": "malware"
    },
    {
      "source_name": "PowMix C2 URL Pattern",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "PowMix",
      "target_type": "malware"
    }
  ]
}