# stix/uac-0247-espionage-campaign-targets-government-and-healthcare-in-ukraine-extraction.json
{
  "_meta": {
    "schema_version": "1.0",
    "generator": "generate-stix-extraction.ts",
    "model": "gemini-2.5-pro",
    "extracted_at": "2026-04-27T17:03:41.467Z",
    "article_slug": "uac-0247-espionage-campaign-targets-government-and-healthcare-in-ukraine",
    "article_date": "2026-04-18",
    "token_usage": {
      "input_tokens": 142025,
      "output_tokens": 4409,
      "cached_tokens": 136033,
      "total_tokens": 146434,
      "cache_hit": true,
      "estimated_cost_usd": 0.068584
    }
  },
  "report_name": "UAC-0247 Espionage Campaign Targets Ukrainian Government and Healthcare with Data-Stealing Malware",
  "threat_actors": [
    {
      "name": "UAC-0247",
      "mitre_id": "",
      "aliases": [],
      "description": "A threat actor group conducting a cyber-espionage campaign against Ukrainian government and healthcare sectors since March 2026, using phishing and a multi-stage infection chain to deploy data-stealing malware.",
      "goals": "Gathering intelligence on Ukrainian government operations and the state of civilian services during a time of conflict.",
      "primary_motivation": "organizational-gain",
      "sophistication": "intermediate"
    }
  ],
  "attack_patterns": [
    {
      "name": "Spearphishing Link",
      "mitre_id": "T1566.002",
      "description": "The attack is initiated via a malicious link in a phishing email, often themed around humanitarian aid."
    },
    {
      "name": "User Execution: Malicious File",
      "mitre_id": "T1204.002",
      "description": "The victim is tricked into opening a malicious .LNK shortcut file contained within a downloaded ZIP archive."
    },
    {
      "name": "System Binary Proxy Execution: Mshta",
      "mitre_id": "T1218.005",
      "description": "The LNK file launches an HTA file, which is executed by mshta.exe to run scripts that can download and execute further payloads."
    },
    {
      "name": "Command and Scripting Interpreter: PowerShell",
      "mitre_id": "T1059.001",
      "description": "The infection chain involves the use of scripting, with PowerShell being a likely interpreter used within the HTA file or subsequent stages."
    },
    {
      "name": "Scheduled Task/Job: Scheduled Task",
      "mitre_id": "T1053.005",
      "description": "The HTA script creates a scheduled task to launch the main executable payload, establishing persistence."
    },
    {
      "name": "Process Injection",
      "mitre_id": "T1055",
      "description": "The malware injects its malicious shellcode into the legitimate Windows process RuntimeBroker.exe to evade detection."
    },
    {
      "name": "Data from Local System",
      "mitre_id": "T1005",
      "description": "The malware payload steals data from the local system, specifically targeting files related to web browsers and the WhatsApp desktop client."
    },
    {
      "name": "Credentials from Password Stores",
      "mitre_id": "T1555.003",
      "description": "The malware specifically targets and steals credential files, cookies, and history from Chromium-based browsers."
    },
    {
      "name": "Exfiltration Over C2 Channel",
      "mitre_id": "T1041",
      "description": "The stolen data is exfiltrated back to an attacker-controlled command-and-control server."
    },
    {
      "name": "Drive-by Compromise",
      "mitre_id": "T1189",
      "description": "Attackers exploited XSS vulnerabilities on legitimate sites to host malicious files, which victims downloaded after clicking a phishing link, constituting a form of drive-by compromise."
    }
  ],
  "malware": [
    {
      "name": "RAVENSHELL-like backdoor",
      "mitre_id": "",
      "type": "backdoor",
      "is_family": false,
      "aliases": [],
      "description": "A data-stealing payload injected into RuntimeBroker.exe, designed to exfiltrate data from Chromium-based browsers and the WhatsApp desktop application. It functions as a reverse TCP shell.",
      "capabilities": [
        "steals-authentication-credentials",
        "exfiltrates-data",
        "communicates-with-c2"
      ]
    },
    {
      "name": "RAVENSHELL",
      "mitre_id": "",
      "type": "backdoor",
      "is_family": true,
      "aliases": [],
      "description": "A known backdoor which the UAC-0247 payload resembles in functionality."
    }
  ],
  "tools": [],
  "vulnerabilities": [
    {
      "name": "Cross-Site Scripting (XSS)",
      "cve": "",
      "description": "Attackers exploited Cross-Site Scripting (XSS) vulnerabilities on legitimate websites to host their malicious ZIP archives for download."
    }
  ],
  "campaigns": [],
  "identities": [
    {
      "name": "Ukrainian Government Organizations",
      "identity_class": "class",
      "sectors": [
        "government"
      ],
      "description": "A primary target sector for the UAC-0247 espionage campaign, including various government bodies in Ukraine."
    },
    {
      "name": "Ukrainian Healthcare Facilities",
      "identity_class": "class",
      "sectors": [
        "healthcare"
      ],
      "description": "A primary target sector for the UAC-0247 espionage campaign, including municipal clinics and hospitals in Ukraine."
    },
    {
      "name": "CERT-UA",
      "identity_class": "organization",
      "sectors": [
        "government",
        "technology"
      ],
      "description": "The Computer Emergency Response Team of Ukraine, which issued the alert about the UAC-0247 campaign."
    }
  ],
  "locations": [
    {
      "name": "Ukraine",
      "country": "Ukraine",
      "region": "eastern-europe"
    }
  ],
  "infrastructure": [
    {
      "name": "UAC-0247 Payload Hosting Infrastructure",
      "infrastructure_types": [
        "hosting-malware"
      ],
      "description": "AI-generated websites or legitimate sites with XSS vulnerabilities used by UAC-0247 to host malicious ZIP archives for download."
    },
    {
      "name": "UAC-0247 C2 Server",
      "infrastructure_types": [
        "command-and-control"
      ],
      "description": "Attacker-controlled server used to receive data exfiltrated by the RAVENSHELL-like backdoor."
    }
  ],
  "observables": [],
  "indicators": [
    {
      "name": "Suspicious HTA File Execution",
      "indicator_types": [
        "malicious-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious HTA File Execution by Email Client or Browser\nid: f1b2c3d4-e5f6-4a7b-8c9d-0e1f2a3b4c5d\nstatus: experimental\ndescription: Detects mshta.exe being launched by a common email client or web browser, which is a TTP used by UAC-0247.\nreferences:\n  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html\nauthor: cyber.netsecops.io\ndate: 2026-04-18\ntags:\n  - attack.t1218.005\n  - attack.execution\nlevel: high\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\mshta.exe'\n    ParentImage|endswith:\n      - '\\outlook.exe'\n      - '\\chrome.exe'\n      - '\\firefox.exe'\n      - '\\msedge.exe'\n  condition: selection",
      "description": "Detects the execution of an HTA file via mshta.exe, a key step in the UAC-0247 infection chain."
    },
    {
      "name": "Suspicious Scheduled Task Creation",
      "indicator_types": [
        "malicious-activity",
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Suspicious Scheduled Task Creation by Script Interpreter\nid: a2b3c4d5-e6f7-4b8c-9d0e-1f2a3b4c5d6e\nstatus: experimental\ndescription: Detects the creation of a scheduled task by a script interpreter like mshta.exe, a persistence technique used by UAC-0247.\nreferences:\n  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html\nauthor: cyber.netsecops.io\ndate: 2026-04-18\ntags:\n  - attack.t1053.005\n  - attack.persistence\nlevel: medium\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\schtasks.exe'\n    CommandLine|contains: '/create'\n    ParentImage|endswith:\n      - '\\mshta.exe'\n      - '\\powershell.exe'\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n  condition: selection",
      "description": "Detects the creation of a scheduled task, which UAC-0247 uses for persistence."
    },
    {
      "name": "Anomalous RuntimeBroker Network Activity",
      "indicator_types": [
        "anomalous-activity"
      ],
      "stix_pattern": "",
      "x_netsecio_sigma_rule": "",
      "pattern": "title: Anomalous RuntimeBroker Network Activity\nid: b3c4d5e6-f7a8-4c9d-0e1f-2a3b4c5d6f7a\nstatus: experimental\ndescription: Detects RuntimeBroker.exe initiating outbound network connections, which is atypical behavior and could indicate process injection as used by UAC-0247.\nreferences:\n  - https://securityaffairs.com/161942/apt/uac-0247-expands-cyber-campaign-across-ukraine.html\nauthor: cyber.netsecops.io\ndate: 2026-04-18\ntags:\n  - attack.t1055\n  - attack.defense_evasion\nlevel: medium\nlogsource:\n  category: network_connection\n  product: windows\ndetection:\n  selection:\n    Image|endswith: '\\RuntimeBroker.exe'\n    Initiated: 'true'\n  filter:\n    DestinationIp|cidr:\n      - '127.0.0.0/8'\n      - '::1/128'\n  condition: selection and not filter",
      "description": "Detects RuntimeBroker.exe making outbound network connections, which is atypical and could indicate process injection as used by UAC-0247."
    }
  ],
  "relationships": [
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "RAVENSHELL-like backdoor",
      "target_type": "malware"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Spearphishing Link",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "User Execution: Malicious File",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "System Binary Proxy Execution: Mshta",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "Drive-by Compromise",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "exploits",
      "target_name": "Cross-Site Scripting (XSS)",
      "target_type": "vulnerability"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukrainian Government Organizations",
      "target_type": "identity"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukrainian Healthcare Facilities",
      "target_type": "identity"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "targets",
      "target_name": "Ukraine",
      "target_type": "location"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "controls",
      "target_name": "UAC-0247 C2 Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "UAC-0247",
      "source_type": "threat-actor",
      "relationship_type": "uses",
      "target_name": "UAC-0247 Payload Hosting Infrastructure",
      "target_type": "infrastructure"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "variant-of",
      "target_name": "RAVENSHELL",
      "target_type": "malware"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Scheduled Task/Job: Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Process Injection",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Data from Local System",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Credentials from Password Stores",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "uses",
      "target_name": "Exfiltration Over C2 Channel",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "RAVENSHELL-like backdoor",
      "source_type": "malware",
      "relationship_type": "beacons-to",
      "target_name": "UAC-0247 C2 Server",
      "target_type": "infrastructure"
    },
    {
      "source_name": "Suspicious HTA File Execution",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "UAC-0247",
      "target_type": "threat-actor"
    },
    {
      "source_name": "Suspicious HTA File Execution",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "System Binary Proxy Execution: Mshta",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Suspicious Scheduled Task Creation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "UAC-0247",
      "target_type": "threat-actor"
    },
    {
      "source_name": "Suspicious Scheduled Task Creation",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Scheduled Task/Job: Scheduled Task",
      "target_type": "attack-pattern"
    },
    {
      "source_name": "Anomalous RuntimeBroker Network Activity",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "UAC-0247",
      "target_type": "threat-actor"
    },
    {
      "source_name": "Anomalous RuntimeBroker Network Activity",
      "source_type": "indicator",
      "relationship_type": "indicates",
      "target_name": "Process Injection",
      "target_type": "attack-pattern"
    }
  ]
}